# peercoin/PLIST
@comment $NetBSD$
sbin/ppcoind
share/examples/peercoin/peercoin.conf
share/examples/rc.d/peercoind
