# boinc-lib/PLIST
@comment $NetBSD: PLIST,v 1.1.1.2 2006/06/18 16:46:05 rumko Exp $
bin/boinc_${PKGVERSION}
bin/boinc_gui
lib/libbenchmark.a
share/examples/rc.d/boinc
