# dnsmasq/PLIST
@comment $NetBSD: PLIST,v 1.1.1.2 2003/10/13 20:09:43 poppnk Exp $
etc/rc.d/dnsmasq
man/man8/dnsmasq.8
sbin/dnsmasq
