# dnsmasq/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2003/07/24 14:16:32 poppnk Exp $
etc/rc.d/dnsmasq
sbin/dnsmasq
share/doc/dnsmasq/doc.html
share/doc/dnsmasq/setup.html
@dirrm share/doc/dnsmasq
