# py-jabber-msnt/PLIST
@comment $NetBSD: PLIST,v 1.8 2006/08/16 15:35:26 schnoebe Exp $
${TRANSPORTDIR}/.svn/entries
${TRANSPORTDIR}/data/defaultJabberAvatar.png
${TRANSPORTDIR}/data/defaultMSNAvatar.png
${TRANSPORTDIR}/src/avatar.py
${TRANSPORTDIR}/src/baseproto/__init__.py
${TRANSPORTDIR}/src/baseproto/glue.py
${TRANSPORTDIR}/src/config.py
${TRANSPORTDIR}/src/contact.py
${TRANSPORTDIR}/src/debug.py
${TRANSPORTDIR}/src/disco.py
${TRANSPORTDIR}/src/ft.py
${TRANSPORTDIR}/src/groupchat.py
${TRANSPORTDIR}/src/housekeep.py
${TRANSPORTDIR}/src/jabw.py
${TRANSPORTDIR}/src/lang.py
${TRANSPORTDIR}/src/legacy/__init__.py
${TRANSPORTDIR}/src/legacy/glue.py
${TRANSPORTDIR}/src/legacy/msn/__init__.py
${TRANSPORTDIR}/src/legacy/msn/msn.py
${TRANSPORTDIR}/src/legacy/msn/msnft.py
${TRANSPORTDIR}/src/legacy/msn/msnp11chl.py
${TRANSPORTDIR}/src/legacy/msn/msnw.py
${TRANSPORTDIR}/src/legacy/msn/test_msn.py
${TRANSPORTDIR}/src/legacy/msn/test_msnw.py
${TRANSPORTDIR}/src/main.py
${TRANSPORTDIR}/src/misciq.py
${TRANSPORTDIR}/src/register.py
${TRANSPORTDIR}/src/session.py
${TRANSPORTDIR}/src/socks5.py
${TRANSPORTDIR}/src/svninfo.py
${TRANSPORTDIR}/src/throttle.py
${TRANSPORTDIR}/src/twistfix/README
${TRANSPORTDIR}/src/twistfix/__init__.py
${TRANSPORTDIR}/src/twistfix/version.py
${TRANSPORTDIR}/src/twistfix/words/__init__.py
${TRANSPORTDIR}/src/twistfix/words/protocols/__init__.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/__init__.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/client.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/component.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/jid.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/jstrports.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/xmpp_stringprep.py
${TRANSPORTDIR}/src/twistfix/words/xish/__init__.py
${TRANSPORTDIR}/src/twistfix/words/xish/domish.py
${TRANSPORTDIR}/src/twistfix/words/xish/utility.py
${TRANSPORTDIR}/src/twistfix/words/xish/xmlstream.py
${TRANSPORTDIR}/src/twistfix/words/xish/xpath.py
${TRANSPORTDIR}/src/twistfix/words/xish/xpathparser.py
${TRANSPORTDIR}/src/utils.py
${TRANSPORTDIR}/src/xdb.py
${TRANSPORTDIR}/src/xmlconfig.py
bin/PyMSNt
share/examples/jabberd/pymsnt.xml
share/examples/rc.d/pymsnt
@dirrm ${TRANSPORTDIR}/src/baseproto
@dirrm ${TRANSPORTDIR}/src/legacy/msn
@dirrm ${TRANSPORTDIR}/src/legacy
@dirrm ${TRANSPORTDIR}/src/twistfix/words/protocols/jabber
@dirrm ${TRANSPORTDIR}/src/twistfix/words/protocols
@dirrm ${TRANSPORTDIR}/src/twistfix/words/xish
@dirrm ${TRANSPORTDIR}/src/twistfix/words
@dirrm ${TRANSPORTDIR}/src/twistfix
@dirrm ${TRANSPORTDIR}/src
@dirrm ${TRANSPORTDIR}/data
@dirrm ${TRANSPORTDIR}/.svn
@dirrm ${TRANSPORTDIR}
