# py-jabber-msnt/PLIST
@comment $NetBSD: PLIST,v 1.5 2006/07/15 00:19:46 schnoebe Exp $
${TRANSPORTDIR}/.svn/entries
${TRANSPORTDIR}/src/avatar.py
${TRANSPORTDIR}/src/baseproto/__init__.py
${TRANSPORTDIR}/src/baseproto/glue.py
${TRANSPORTDIR}/src/config.py
${TRANSPORTDIR}/src/contact.py
${TRANSPORTDIR}/src/debug.py
${TRANSPORTDIR}/src/disco.py
${TRANSPORTDIR}/src/ft.py
${TRANSPORTDIR}/src/groupchat.py
${TRANSPORTDIR}/src/housekeep.py
${TRANSPORTDIR}/src/jabw.py
${TRANSPORTDIR}/src/lang.py
${TRANSPORTDIR}/src/legacy/__init__.py
${TRANSPORTDIR}/src/legacy/glue.py
${TRANSPORTDIR}/src/main.py
${TRANSPORTDIR}/src/misciq.py
${TRANSPORTDIR}/src/register.py
${TRANSPORTDIR}/src/session.py
${TRANSPORTDIR}/src/tlib/__init__.py
${TRANSPORTDIR}/src/tlib/domish.py
${TRANSPORTDIR}/src/tlib/jabber/__init__.py
${TRANSPORTDIR}/src/tlib/jabber/client.py
${TRANSPORTDIR}/src/tlib/jabber/component.py
${TRANSPORTDIR}/src/tlib/jabber/jid.py
${TRANSPORTDIR}/src/tlib/jabber/jstrports.py
${TRANSPORTDIR}/src/tlib/jabber/xmpp_stringprep.py
${TRANSPORTDIR}/src/tlib/msn/__init__.py
${TRANSPORTDIR}/src/tlib/msn/msn.py
${TRANSPORTDIR}/src/tlib/msn/msnft.py
${TRANSPORTDIR}/src/tlib/msn/msnp11chl.py
${TRANSPORTDIR}/src/tlib/msn/msnw.py
${TRANSPORTDIR}/src/tlib/msn/test_msn.py
${TRANSPORTDIR}/src/tlib/msn/test_msnw.py
${TRANSPORTDIR}/src/tlib/socks5.py
${TRANSPORTDIR}/src/tlib/throttle.py
${TRANSPORTDIR}/src/tlib/xmlstream.py
${TRANSPORTDIR}/src/tlib/xmlw.py
${TRANSPORTDIR}/src/utils.py
${TRANSPORTDIR}/src/xdb.py
${TRANSPORTDIR}/src/xmlconfig.py
${TRANSPORTDIR}/data/defaultJabberAvatar.png
${TRANSPORTDIR}/data/defaultMSNAvatar.png
bin/PyMSNt
share/examples/jabberd/pymsnt.xml
share/examples/rc.d/pymsnt
@dirrm ${TRANSPORTDIR}/src/baseproto
@dirrm ${TRANSPORTDIR}/src/legacy
@dirrm ${TRANSPORTDIR}/src/tlib/jabber
@dirrm ${TRANSPORTDIR}/src/tlib/msn
@dirrm ${TRANSPORTDIR}/src/tlib
@dirrm ${TRANSPORTDIR}/src
@dirrm ${TRANSPORTDIR}/data
@dirrm ${TRANSPORTDIR}
