# libreswan-git/PLIST
@comment $NetBSD$
libexec/ipsec/_asn1check
libexec/ipsec/_dncheck
libexec/ipsec/_enumcheck
libexec/ipsec/_hunkcheck
libexec/ipsec/_import_crl
libexec/ipsec/_ipcheck
libexec/ipsec/_jambufcheck
libexec/ipsec/_keyidcheck
libexec/ipsec/_plutorun
libexec/ipsec/_secretcensor
libexec/ipsec/_timecheck
libexec/ipsec/_ttodatacheck
libexec/ipsec/_unbound-hook
libexec/ipsec/_updown
libexec/ipsec/_updown.bsd
libexec/ipsec/_vendoridcheck
libexec/ipsec/addconn
libexec/ipsec/algparse
libexec/ipsec/barf
libexec/ipsec/cavp
libexec/ipsec/ecdsasigkey
libexec/ipsec/letsencrypt
libexec/ipsec/newhostkey
libexec/ipsec/pluto
libexec/ipsec/readwriteconf
libexec/ipsec/rsasigkey
libexec/ipsec/setup
libexec/ipsec/showhostkey
libexec/ipsec/showroute
libexec/ipsec/whack
man/man5/ipsec.conf.5
man/man5/ipsec.secrets.5
man/man7/ipsec_pluto.7
man/man8/ipsec.8
man/man8/ipsec__import_crl.8
man/man8/ipsec__plutorun.8
man/man8/ipsec__secretcensor.8
man/man8/ipsec__unbound-hook.8
man/man8/ipsec__updown.8
man/man8/ipsec__updown.bsdkame.8
man/man8/ipsec_add.8
man/man8/ipsec_addconn.8
man/man8/ipsec_barf.8
man/man8/ipsec_briefstatus.8
man/man8/ipsec_certutil.8
man/man8/ipsec_checkconfig.8
man/man8/ipsec_checknflog.8
man/man8/ipsec_checknss.8
man/man8/ipsec_connectionstatus.8
man/man8/ipsec_crlutil.8
man/man8/ipsec_delete.8
man/man8/ipsec_down.8
man/man8/ipsec_ecdsasigkey.8
man/man8/ipsec_fetchcrls.8
man/man8/ipsec_fipsstatus.8
man/man8/ipsec_globalstatus.8
man/man8/ipsec_import.8
man/man8/ipsec_initnss.8
man/man8/ipsec_letsencrypt.8
man/man8/ipsec_listall.8
man/man8/ipsec_listcacerts.8
man/man8/ipsec_listcrls.8
man/man8/ipsec_listen.8
man/man8/ipsec_listpubkeys.8
man/man8/ipsec_modutil.8
man/man8/ipsec_newhostkey.8
man/man8/ipsec_ondemand.8
man/man8/ipsec_pk12util.8
man/man8/ipsec_pluto.8
man/man8/ipsec_purgeocsp.8
man/man8/ipsec_readwriteconf.8
man/man8/ipsec_redirect.8
man/man8/ipsec_replace.8
man/man8/ipsec_rereadall.8
man/man8/ipsec_rereadcerts.8
man/man8/ipsec_rereadsecrets.8
man/man8/ipsec_restart.8
man/man8/ipsec_route.8
man/man8/ipsec_rsasigkey.8
man/man8/ipsec_setup.8
man/man8/ipsec_showhostkey.8
man/man8/ipsec_showroute.8
man/man8/ipsec_showstates.8
man/man8/ipsec_shuntstatus.8
man/man8/ipsec_start.8
man/man8/ipsec_status.8
man/man8/ipsec_stop.8
man/man8/ipsec_trafficstatus.8
man/man8/ipsec_unroute.8
man/man8/ipsec_up.8
man/man8/ipsec_vfychain.8
man/man8/ipsec_whack.8
man/man8/pluto.8
sbin/ipsec
share/examples/libreswan/ipsec.conf-sample
share/examples/libreswan/ipsec.d/policies/block
share/examples/libreswan/ipsec.d/policies/clear
share/examples/libreswan/ipsec.d/policies/clear-or-private
share/examples/libreswan/ipsec.d/policies/portexcludes.conf
share/examples/libreswan/ipsec.d/policies/private
share/examples/libreswan/ipsec.d/policies/private-or-clear
share/examples/libreswan/ipsec.secrets-sample
share/examples/libreswan/logrotate.d/libreswan
share/examples/libreswan/pam.d/pluto
share/examples/libreswan/rc.d/pluto.sh
