# py-jabber-aimt/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/10/11 10:45:10 thomasklausner Exp $
bin/PyAIMt
${TRANSPORTDIR}/src/avatar.py
${TRANSPORTDIR}/src/adhoc.py
${TRANSPORTDIR}/src/config.py
${TRANSPORTDIR}/src/bool.py
${TRANSPORTDIR}/src/langs/__init__.py
${TRANSPORTDIR}/src/langs/cs.py
${TRANSPORTDIR}/src/langs/de.py
${TRANSPORTDIR}/src/langs/en.py
${TRANSPORTDIR}/src/langs/fr.py
${TRANSPORTDIR}/src/langs/nl.py
${TRANSPORTDIR}/src/langs/pl.py
${TRANSPORTDIR}/src/langs/ru.py
${TRANSPORTDIR}/src/utils.py
${TRANSPORTDIR}/src/contact.py
${TRANSPORTDIR}/src/debug.py
${TRANSPORTDIR}/src/disco.py
${TRANSPORTDIR}/src/globals.py
${TRANSPORTDIR}/src/groupchat.py
${TRANSPORTDIR}/src/imgmanip.py
${TRANSPORTDIR}/src/iq.py
${TRANSPORTDIR}/src/jabw.py
${TRANSPORTDIR}/src/lang.py
${TRANSPORTDIR}/src/legacy/services/AIMURITranslate.py
${TRANSPORTDIR}/src/legacy/services/ChangeEmail.py
${TRANSPORTDIR}/src/legacy/services/ChangePassword.py
${TRANSPORTDIR}/src/legacy/services/ConfirmAccount.py
${TRANSPORTDIR}/src/legacy/services/EmailLookup.py
${TRANSPORTDIR}/src/legacy/services/FormatScreenName.py
${TRANSPORTDIR}/src/legacy/services/__init__.py
${TRANSPORTDIR}/src/legacy/__init__.py
${TRANSPORTDIR}/src/legacy/aimt.py
${TRANSPORTDIR}/src/legacy/buddies.py
${TRANSPORTDIR}/src/legacy/glue.py
${TRANSPORTDIR}/src/register.py
${TRANSPORTDIR}/src/main.py
${TRANSPORTDIR}/src/services/GatewayTranslator.py
${TRANSPORTDIR}/src/services/ConnectUsers.py
${TRANSPORTDIR}/src/services/EntityTime.py
${TRANSPORTDIR}/src/services/IqAvatarFactory.py
${TRANSPORTDIR}/src/services/LastActivity.py
${TRANSPORTDIR}/src/services/PingService.py
${TRANSPORTDIR}/src/services/RosterRetrieval.py
${TRANSPORTDIR}/src/services/SearchFactory.py
${TRANSPORTDIR}/src/services/Statistics.py
${TRANSPORTDIR}/src/services/VCardFactory.py
${TRANSPORTDIR}/src/services/VersionTeller.py
${TRANSPORTDIR}/src/services/__init__.py
${TRANSPORTDIR}/src/sasl.py
${TRANSPORTDIR}/src/tlib/countrycodes.py
${TRANSPORTDIR}/src/tlib/__init__.py
${TRANSPORTDIR}/src/tlib/httpcompat.py
${TRANSPORTDIR}/src/tlib/ifcompat.py
${TRANSPORTDIR}/src/tlib/oscar.py
${TRANSPORTDIR}/src/tlib/scheduler.py
${TRANSPORTDIR}/src/tlib/socks5.py
${TRANSPORTDIR}/src/tlib/sockserror.py
${TRANSPORTDIR}/src/web/__init__.py
${TRANSPORTDIR}/src/web/handler.py
${TRANSPORTDIR}/src/web/xmppcred.py
${TRANSPORTDIR}/src/session.py
${TRANSPORTDIR}/src/svninfo.py
${TRANSPORTDIR}/src/xdb/__init__.py
${TRANSPORTDIR}/src/xdb/mysql.py
${TRANSPORTDIR}/src/xdb/template.py
${TRANSPORTDIR}/src/xdb/legacyaimtransport.py
${TRANSPORTDIR}/src/xdb/xmlfiles.py
${TRANSPORTDIR}/src/twistfix/__init__.py
${TRANSPORTDIR}/src/twistfix/README
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/xmpp_stringprep.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/__init__.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/client.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/component.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/jid.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/jstrports.py
${TRANSPORTDIR}/src/twistfix/words/protocols/__init__.py
${TRANSPORTDIR}/src/twistfix/words/__init__.py
${TRANSPORTDIR}/src/twistfix/words/xish/xpathparser.py
${TRANSPORTDIR}/src/twistfix/words/xish/__init__.py
${TRANSPORTDIR}/src/twistfix/words/xish/domish.py
${TRANSPORTDIR}/src/twistfix/words/xish/utility.py
${TRANSPORTDIR}/src/twistfix/words/xish/xmlstream.py
${TRANSPORTDIR}/src/twistfix/words/xish/xpath.py
${TRANSPORTDIR}/src/twistfix/version.py
${TRANSPORTDIR}/src/xmlconfig.py
${TRANSPORTDIR}/data/www/css/style.css
${TRANSPORTDIR}/data/www/images/pyaim.png
${TRANSPORTDIR}/data/www/template.html
${TRANSPORTDIR}/data/defaultAIMAvatar.png
${TRANSPORTDIR}/data/defaultICQAvatar.png
${TRANSPORTDIR}/tools/db-setup.mysql
${TRANSPORTDIR}/tools/Makefile
${TRANSPORTDIR}/tools/infodump.py
${TRANSPORTDIR}/tools/managessi.py
${TRANSPORTDIR}/tools/migrate.py
${JABBER_EXDIR}/pyaimt.xml
share/examples/rc.d/pyaimt
