# bcmfw/PLIST
@comment $NetBSD: PLIST,v 1.4 2006/08/28 11:48:26 plunky Exp $
man/man8/bcmfw.8
sbin/bcmfw
share/examples/rc.d/bcmfw
${FIRMWARE_DIR}/BCM-LEGAL.txt
${FIRMWARE_DIR}/BCM2033-FW.bin
${FIRMWARE_DIR}/BCM2033-MD.hex
@dirrm ${FIRMWARE_DIR}
