# dcc/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2008/12/22 17:26:14 pebog Exp $
bin/cdcc
bin/dccproc
libexec/cron-dccd
${PLIST.dccd}libexec/dbclean
${PLIST.dccd}libexec/dblist
libexec/dcc-stats-collect
libexec/dcc-stats-graph
libexec/dcc-stats-init
${PLIST.dccd}libexec/dccd
${PLIST.dccifd}libexec/dccifd
${PLIST.dccm}libexec/dccm
libexec/dccsight
${PLIST.dns-helper}libexec/dns-helper
libexec/fetch-testmsg-whitelist
libexec/fetchblack
libexec/fetchids
libexec/hackmc
libexec/list-clients
libexec/logger
libexec/newwebuser
libexec/rcDCC
libexec/start-dccd
libexec/start-dccifd
libexec/start-dccm
libexec/start-grey
libexec/stats-get
libexec/stop-dccd
libexec/uninstalldcc
libexec/updatedcc
libexec/wlist
man/cat8/cdcc.0
man/cat8/dbclean.0
man/cat8/dblist.0
man/cat8/dcc.0
man/cat8/dccd.0
man/cat8/dccifd.0
man/cat8/dccm.0
man/cat8/dccproc.0
man/cat8/dccsight.0
share/examples/dcc/cgi-bin/README
share/examples/dcc/cgi-bin/chgpasswd
share/examples/dcc/cgi-bin/common
share/examples/dcc/cgi-bin/edit-whiteclnt
share/examples/dcc/cgi-bin/footer
share/examples/dcc/cgi-bin/header
share/examples/dcc/cgi-bin/http2https
share/examples/dcc/cgi-bin/list-log
share/examples/dcc/cgi-bin/list-msg
share/examples/dcc/cgi-bin/webuser-notify
share/examples/dcc/dcc_conf
share/examples/dcc/flod
share/examples/dcc/grey_flod
share/examples/dcc/grey_whitelist
share/examples/dcc/ids
share/examples/dcc/map.txt
share/examples/dcc/whiteclnt
share/examples/dcc/whitecommon
share/examples/dcc/whitelist
${PLIST.dccd}share/examples/rc.d/dccd
${PLIST.dccifd}share/examples/rc.d/dccifd
${PLIST.dccm}share/examples/rc.d/dccm
@dirrm share/examples/dcc/cgi-bin
@dirrm share/examples/dcc
