# ypbind-mt/PLIST
@comment $NetBSD: PLIST,v 1.3 2006/04/18 00:46:39 jlamwww Exp $
man/man5/yp.conf.5
man/man8/ypbind.8
sbin/ypbind
share/examples/rc.d/ypbind
share/locale/de/LC_MESSAGES/ypbind-mt.mo
