# ypbind-mt/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:08 jeremy-c-reed Exp $
man/man5/yp.conf.5
man/man8/ypbind.8
sbin/ypbind
${PKGLOCALEDIR}/locale/de/LC_MESSAGES/ypbind-mt.mo
share/examples/rc.d/ypbind
