# ushare/PLIST
@comment $NetBSD: PLIST,v 1.6 2011/01/20 08:27:11 obache Exp $
bin/ushare
man/man1/ushare.1
share/examples/rc.d/ushare
share/examples/ushare/init.d/ushare
share/examples/ushare/ushare.conf
share/locale/de/LC_MESSAGES/ushare.mo
share/locale/fr/LC_MESSAGES/ushare.mo
