# py-jabber-palaver/PLIST
@comment $NetBSD: PLIST,v 1.7 2009/10/11 10:45:10 thomasklausner Exp $
${PYSITELIB}/palaver/__init__.py
${PYSITELIB}/palaver/__init__.pyc
${PYSITELIB}/palaver/__init__.pyo
${PYSITELIB}/palaver/admin.py
${PYSITELIB}/palaver/admin.pyc
${PYSITELIB}/palaver/admin.pyo
${PYSITELIB}/palaver/cred.py
${PYSITELIB}/palaver/cred.pyc
${PYSITELIB}/palaver/cred.pyo
${PYSITELIB}/palaver/dir_storage.py
${PYSITELIB}/palaver/dir_storage.pyc
${PYSITELIB}/palaver/dir_storage.pyo
${PYSITELIB}/palaver/groupchat.py
${PYSITELIB}/palaver/groupchat.pyc
${PYSITELIB}/palaver/groupchat.pyo
${PYSITELIB}/palaver/memory_storage.py
${PYSITELIB}/palaver/memory_storage.pyc
${PYSITELIB}/palaver/memory_storage.pyo
${PYSITELIB}/palaver/muc.py
${PYSITELIB}/palaver/muc.pyc
${PYSITELIB}/palaver/muc.pyo
${PYSITELIB}/palaver/palaver.py
${PYSITELIB}/palaver/palaver.pyc
${PYSITELIB}/palaver/palaver.pyo
${PYSITELIB}/palaver/pgsql_storage.py
${PYSITELIB}/palaver/pgsql_storage.pyc
${PYSITELIB}/palaver/pgsql_storage.pyo
${PYSITELIB}/palaver/plog.py
${PYSITELIB}/palaver/plog.pyc
${PYSITELIB}/palaver/plog.pyo
${PYSITELIB}/palaver/storage.py
${PYSITELIB}/palaver/storage.pyc
${PYSITELIB}/palaver/storage.pyo
${PYSITELIB}/palaver/tap.py
${PYSITELIB}/palaver/tap.pyc
${PYSITELIB}/palaver/tap.pyo
${PYSITELIB}/palaver/test/__init__.py
${PYSITELIB}/palaver/test/__init__.pyc
${PYSITELIB}/palaver/test/__init__.pyo
${PYSITELIB}/palaver/test/harness.py
${PYSITELIB}/palaver/test/harness.pyc
${PYSITELIB}/palaver/test/harness.pyo
${PYSITELIB}/palaver/test/pgsql_storage.py
${PYSITELIB}/palaver/test/pgsql_storage.pyc
${PYSITELIB}/palaver/test/pgsql_storage.pyo
${PYSITELIB}/palaver/test/readlog.py
${PYSITELIB}/palaver/test/readlog.pyc
${PYSITELIB}/palaver/test/readlog.pyo
${PYSITELIB}/palaver/test/regression.py
${PYSITELIB}/palaver/test/regression.pyc
${PYSITELIB}/palaver/test/regression.pyo
${PYSITELIB}/palaver/test/storage.py
${PYSITELIB}/palaver/test/storage.pyc
${PYSITELIB}/palaver/test/storage.pyo
${PYSITELIB}/palaver/test/xep045.py
${PYSITELIB}/palaver/test/xep045.pyc
${PYSITELIB}/palaver/test/xep045.pyo
${PYSITELIB}/palaver/xmpp/__init__.py
${PYSITELIB}/palaver/xmpp/__init__.pyc
${PYSITELIB}/palaver/xmpp/__init__.pyo
${PYSITELIB}/palaver/xmpp/disco.py
${PYSITELIB}/palaver/xmpp/disco.pyc
${PYSITELIB}/palaver/xmpp/disco.pyo
${PYSITELIB}/palaver/xmpp/error.py
${PYSITELIB}/palaver/xmpp/error.pyc
${PYSITELIB}/palaver/xmpp/error.pyo
${PYSITELIB}/palaver/xmpp/ns.py
${PYSITELIB}/palaver/xmpp/ns.pyc
${PYSITELIB}/palaver/xmpp/ns.pyo
${PYSITELIB}/twisted/plugins/palaver.py
${PYSITELIB}/twisted/plugins/palaver.pyc
${PYSITELIB}/twisted/plugins/palaver.pyo
share/doc/palaver/README.txt
share/doc/palaver/XEP.txt
share/doc/palaver/muc-pgsql.sql
share/examples/jabber/palaver.xml
share/examples/rc.d/palaver
@exec "rmdir share/examples/jabber || true"
