# hal/PLIST
@comment $NetBSD: PLIST,v 1.2 2007/12/02 21:19:23 ondratu Exp $
bin/hal-device
bin/hal-disable-polling
bin/hal-find-by-capability
bin/hal-find-by-property
bin/hal-get-property
bin/hal-is-caller-locked-out
bin/hal-lock
bin/hal-set-property
bin/lshal
include/hal/libhal.h
include/hal/libhal-storage.h
@comment this addons could not be build, so this could be generate dynamicle
libexec/hald-addon-acpi
libexec/hald-addon-cpufreq
libexec/hald-addon-hid-ups
libexec/hald-addon-input
libexec/hald-addon-macbook-backlight
libexec/hald-addon-macbookpro-backlight
libexec/hald-addon-pmu
libexec/hald-addon-storage
libexec/hald-addon-usb-csr
libexec/hald-generate-fdi-cache
libexec/hald-probe-hiddev
libexec/hald-probe-ieee1394-unit
libexec/hald-probe-input
libexec/hald-probe-net-bluetooth
libexec/hald-probe-pc-floppy
libexec/hald-probe-printer
libexec/hald-probe-serial
libexec/hald-probe-smbios
libexec/hald-probe-storage
libexec/hald-probe-volume
libexec/hald-runner
libexec/hal-ipw-killswitch-${LOWER_OPSYS}
libexec/hal-storage-cleanup-all-mountpoints
libexec/hal-storage-cleanup-mountpoint
libexec/hal-storage-closetray
libexec/hal-storage-eject
libexec/hal-storage-mount
libexec/hal-storage-unmount
libexec/hal-system-power-pm-is-supported
libexec/hal-system-power-pmu
libexec/hal-system-sonypic
lib/hal/scripts/hal-functions
lib/hal/scripts/hal-luks-remove
lib/hal/scripts/hal-luks-setup
lib/hal/scripts/hal-luks-teardown
lib/hal/scripts/hal-system-killswitch-get-power
lib/hal/scripts/hal-system-killswitch-set-power
lib/hal/scripts/hal-system-lcd-get-brightness
lib/hal/scripts/hal-system-lcd-set-brightness
lib/hal/scripts/hal-system-power-hibernate
lib/hal/scripts/hal-system-power-reboot
lib/hal/scripts/hal-system-power-set-power-save
lib/hal/scripts/hal-system-power-shutdown
lib/hal/scripts/hal-system-power-suspend
lib/hal/scripts/hal-system-power-suspend-hybrid
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-luks-remove-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-luks-setup-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-luks-teardown-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-killswitch-get-power-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-killswitch-set-power-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-lcd-get-brightness-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-lcd-set-brightness-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-power-hibernate-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-power-reboot-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-power-set-power-save-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-power-shutdown-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-power-suspend-hybrid-${LOWER_OPSYS}
${SS}lib/hal/scripts/${LOWER_OPSYS}/hal-system-power-suspend-${LOWER_OPSYS}
lib/libhal.la
lib/libhal-storage.la
lib/pkgconfig/hal.pc
lib/pkgconfig/hal-storage.pc
man/man1/hal-disable-polling.1
man/man1/hal-find-by-capability.1
man/man1/hal-find-by-property.1
man/man1/hal-get-property.1
man/man1/hal-is-caller-locked-out.1
man/man1/hal-lock.1
man/man1/hal-set-property.1
man/man1/lshal.1
man/man8/hald.8
sbin/hald
share/examples/hal/hal.conf
share/examples/hal/90-hal.rules
share/examples/rc.d/hald
share/gtk-doc/html/libhal/home.png
share/gtk-doc/html/libhal/index.html
share/gtk-doc/html/libhal/index.sgml
share/gtk-doc/html/libhal/ix01.html
share/gtk-doc/html/libhal/left.png
share/gtk-doc/html/libhal/libhal.devhelp
share/gtk-doc/html/libhal/libhal.devhelp2
share/gtk-doc/html/libhal/libhal-libhal.html
share/gtk-doc/html/libhal/license.html
share/gtk-doc/html/libhal/right.png
share/gtk-doc/html/libhal/rn01.html
share/gtk-doc/html/libhal-storage/home.png
share/gtk-doc/html/libhal-storage/index.html
share/gtk-doc/html/libhal-storage/index.sgml
share/gtk-doc/html/libhal-storage/ix01.html
share/gtk-doc/html/libhal-storage/left.png
share/gtk-doc/html/libhal-storage/libhal-storage.devhelp
share/gtk-doc/html/libhal-storage/libhal-storage.devhelp2
share/gtk-doc/html/libhal-storage/libhal-storage-libhal-storage.html
share/gtk-doc/html/libhal-storage/license.html
share/gtk-doc/html/libhal-storage/right.png
share/gtk-doc/html/libhal-storage/rn01.html
share/gtk-doc/html/libhal-storage/style.css
share/gtk-doc/html/libhal-storage/up.png
share/gtk-doc/html/libhal-storage/version.xml
share/gtk-doc/html/libhal/style.css
share/gtk-doc/html/libhal/up.png
share/hal/fdi/fdi.dtd
share/hal/fdi/policy/10osvendor/10-cpufreq.fdi
share/hal/fdi/policy/10osvendor/10-input-policy.fdi
share/hal/fdi/policy/10osvendor/10-laptop-panel-mgmt-policy.fdi
share/hal/fdi/policy/10osvendor/10-macbook-backlight.fdi
share/hal/fdi/policy/10osvendor/10-macbookpro-utils.fdi
share/hal/fdi/policy/10osvendor/10-power-mgmt-policy.fdi
share/hal/fdi/policy/10osvendor/10-rfkill-switch.fdi
share/hal/fdi/policy/10osvendor/10-usbcsr-mice.fdi
share/hal/fdi/policy/10osvendor/10-x11-input.fdi
share/hal/fdi/policy/10osvendor/15-storage-luks.fdi
share/hal/fdi/policy/10osvendor/20-storage-methods.fdi
@exec ${MKDIR} %D/share/hal/fdi/preprobe/20thirdparty
@dirrm share/hal/fdi/preprobe/20thirdparty
@exec ${MKDIR} %D/share/hal/fdi/preprobe/10osvendor
@dirrm share/hal/fdi/preprobe/10osvendor
@dirrm share/hal/fdi/preprobe
@exec ${MKDIR} %D/share/hal/fdi/policy/20thirdparty
@dirrm share/hal/fdi/policy/20thirdparty
@dirrm share/hal/fdi/policy/10osvendor
@dirrm share/hal/fdi/policy
@exec ${MKDIR} %D/share/hal/fdi/information/20thirdparty
@dirrm share/hal/fdi/information/20thirdparty
@exec ${MKDIR} %D/share/hal/fdi/information/10freedesktop
@dirrm share/hal/fdi/information/10freedesktop
@dirrm share/hal/fdi/information
@dirrm share/hal/fdi
@dirrm share/hal
@dirrm share/gtk-doc/html/libhal-storage
@dirrm share/gtk-doc/html/libhal
@dirrm share/examples/hal
@dirrm lib/hal/scripts/${LOWER_OPSYS}
@dirrm lib/hal/scripts
@dirrm lib/hal
@dirrm include/hal
@exec ${MKDIR} %D/etc/hal/fdi/preprobe
@dirrm etc/hal/fdi/preprobe
@exec ${MKDIR} %D/etc/hal/fdi/policy
@dirrm etc/hal/fdi/policy
@exec ${MKDIR} %D/etc/hal/fdi/information
@dirrm etc/hal/fdi/information
@dirrm etc/hal/fdi
@dirrm etc/hal
