# haproxy15/PLIST
@comment $NetBSD$
man/man1/haproxy.1
sbin/haproxy
sbin/haproxy-systemd-wrapper
share/doc/haproxy/architecture.txt
share/doc/haproxy/configuration.txt
share/doc/haproxy/haproxy-en.txt
share/doc/haproxy/haproxy-fr.txt
share/examples/haproxy/haproxy.cfg
${PLIST.rcd}share/examples/rc.d/haproxy
