# wicd/PLIST
@comment $NetBSD: PLIST,v 1.4 2012/05/10 20:28:34 noud4 Exp $
bin/wicd-client
bin/wicd-curses
${PYSITELIB}/wicd/backend.pyc
${PYSITELIB}/wicd/backend.py
${PYSITELIB}/wicd/backend.pyo
${PYSITELIB}/wicd/configmanager.pyc
${PYSITELIB}/wicd/configmanager.py
${PYSITELIB}/wicd/configmanager.pyo
${PYSITELIB}/wicd/dbusmanager.pyc
${PYSITELIB}/wicd/dbusmanager.py
${PYSITELIB}/wicd/dbusmanager.pyo
${PYSITELIB}/wicd/gui.pyc
${PYSITELIB}/wicd/gui.py
${PYSITELIB}/wicd/gui.pyo
${PYSITELIB}/wicd/guiutil.pyc
${PYSITELIB}/wicd/guiutil.py
${PYSITELIB}/wicd/guiutil.pyo
${PYSITELIB}/wicd/__init__.pyc
${PYSITELIB}/wicd/__init__.py
${PYSITELIB}/wicd/__init__.pyo
${PYSITELIB}/wicd/logfile.pyc
${PYSITELIB}/wicd/logfile.py
${PYSITELIB}/wicd/logfile.pyo
${PYSITELIB}/wicd/misc.pyc
${PYSITELIB}/wicd/misc.py
${PYSITELIB}/wicd/misc.pyo
${PYSITELIB}/wicd/netentry.pyc
${PYSITELIB}/wicd/netentry.py
${PYSITELIB}/wicd/netentry.pyo
${PYSITELIB}/wicd/networking.pyc
${PYSITELIB}/wicd/networking.py
${PYSITELIB}/wicd/networking.pyo
${PYSITELIB}/wicd/prefs.pyc
${PYSITELIB}/wicd/prefs.py
${PYSITELIB}/wicd/prefs.pyo
${PYSITELIB}/wicd/translations.pyc
${PYSITELIB}/wicd/translations.py
${PYSITELIB}/wicd/translations.pyo
${PYSITELIB}/wicd/wnettools.pyc
${PYSITELIB}/wicd/wnettools.py
${PYSITELIB}/wicd/wnettools.pyo
${PYSITELIB}/wicd/wpath.py
${PYSITELIB}/wicd/wpath.pyc
${PYSITELIB}/wicd/wpath.pyo
${PYSITELIB}/${EGG_INFODIR}
lib/wicd/curses_misc.py
lib/wicd/autoconnect.py
lib/wicd/configscript_curses.py
lib/wicd/configscript.py
lib/wicd/monitor.py
lib/wicd/netentry_curses.py
lib/wicd/prefs_curses.py
lib/wicd/suspend.py
lib/wicd/wicd-client.py
lib/wicd/wicd-curses.py
lib/wicd/wicd-daemon.py
lib/wicd/backends/be-external.py
lib/wicd/backends/be-ioctl.py
man/man5/wicd-manager-settings.conf.5
man/man5/wicd-wired-settings.conf.5
man/man5/wicd-wireless-settings.conf.5
man/man8/wicd.8
man/man8/wicd-curses.8
sbin/wicd
share/locale/sl/LC_MESSAGES/wicd.mo
share/locale/es_MX/LC_MESSAGES/wicd.mo
share/locale/te/LC_MESSAGES/wicd.mo
share/locale/id/LC_MESSAGES/wicd.mo
share/locale/de/LC_MESSAGES/wicd.mo
share/locale/fa/LC_MESSAGES/wicd.mo
share/locale/sk/LC_MESSAGES/wicd.mo
share/locale/uk/LC_MESSAGES/wicd.mo
share/locale/es/LC_MESSAGES/wicd.mo
share/locale/vi/LC_MESSAGES/wicd.mo
share/locale/ml/LC_MESSAGES/wicd.mo
share/locale/ko/LC_MESSAGES/wicd.mo
share/locale/sr/LC_MESSAGES/wicd.mo
share/locale/ro/LC_MESSAGES/wicd.mo
share/locale/nl/LC_MESSAGES/wicd.mo
share/locale/de_DE/LC_MESSAGES/wicd.mo
share/locale/el/LC_MESSAGES/wicd.mo
share/locale/fr_CA/LC_MESSAGES/wicd.mo
share/locale/zh_CN/LC_MESSAGES/wicd.mo
share/locale/it/LC_MESSAGES/wicd.mo
share/locale/fr/LC_MESSAGES/wicd.mo
share/locale/pl/LC_MESSAGES/wicd.mo
share/locale/zh_TW/LC_MESSAGES/wicd.mo
share/locale/et/LC_MESSAGES/wicd.mo
share/locale/es_AR/LC_MESSAGES/wicd.mo
share/locale/es_ES/LC_MESSAGES/wicd.mo
share/locale/ka/LC_MESSAGES/wicd.mo
share/locale/pt/LC_MESSAGES/wicd.mo
share/locale/gl/LC_MESSAGES/wicd.mo
share/locale/bg/LC_MESSAGES/wicd.mo
share/locale/es_GT/LC_MESSAGES/wicd.mo
share/locale/pt_BR/LC_MESSAGES/wicd.mo
share/locale/fi/LC_MESSAGES/wicd.mo
share/locale/eo/LC_MESSAGES/wicd.mo
share/locale/es_NI/LC_MESSAGES/wicd.mo
share/locale/nl_NL/LC_MESSAGES/wicd.mo
share/locale/eu/LC_MESSAGES/wicd.mo
share/locale/zh_HK/LC_MESSAGES/wicd.mo
share/locale/es_VE/LC_MESSAGES/wicd.mo
share/locale/no/LC_MESSAGES/wicd.mo
share/locale/ar_EG/LC_MESSAGES/wicd.mo
share/locale/lv/LC_MESSAGES/wicd.mo
share/locale/tr/LC_MESSAGES/wicd.mo
share/locale/es_CL/LC_MESSAGES/wicd.mo
share/locale/cs/LC_MESSAGES/wicd.mo
share/locale/hu/LC_MESSAGES/wicd.mo
share/locale/ru_RU/LC_MESSAGES/wicd.mo
share/locale/da/LC_MESSAGES/wicd.mo
share/locale/sv/LC_MESSAGES/wicd.mo
share/locale/lt/LC_MESSAGES/wicd.mo
share/locale/it_IT/LC_MESSAGES/wicd.mo
share/locale/ru/LC_MESSAGES/wicd.mo
share/locale/ja/LC_MESSAGES/wicd.mo
share/locale/he/LC_MESSAGES/wicd.mo
share/locale/kk/LC_MESSAGES/wicd.mo
share/locale/ca/LC_MESSAGES/wicd.mo
share/icons/hicolor/24x24/apps/wicd-client.png
share/icons/hicolor/96x96/apps/wicd-client.png
share/icons/hicolor/48x48/apps/wicd-client.png
share/icons/hicolor/128x128/apps/wicd-client.png
share/icons/hicolor/22x22/apps/wicd-client.png
share/icons/hicolor/scalable/apps/wicd-client.svg
share/icons/hicolor/32x32/apps/wicd-client.png
share/icons/hicolor/64x64/apps/wicd-client.png
share/icons/hicolor/72x72/apps/wicd-client.png
share/icons/hicolor/16x16/apps/wicd-client.png
share/icons/hicolor/36x36/apps/wicd-client.png
share/icons/hicolor/192x192/apps/wicd-client.png
share/wicd/wicd.glade
man/man1/wicd-client.1
share/doc/wicd/AUTHORS
share/doc/wicd/CHANGES
share/doc/wicd/INSTALL
share/doc/wicd/LICENSE
share/doc/wicd/README
share/pixmaps/wicd/bad-signal-lock.png
share/pixmaps/wicd/bad-signal.png
share/pixmaps/wicd/both-bad-signal-lock.png
share/pixmaps/wicd/both-bad-signal.png
share/pixmaps/wicd/both-good-signal-lock.png
share/pixmaps/wicd/both-good-signal.png
share/pixmaps/wicd/both-high-signal-lock.png
share/pixmaps/wicd/both-high-signal.png
share/pixmaps/wicd/both-low-signal-lock.png
share/pixmaps/wicd/both-low-signal.png
share/pixmaps/wicd/good-signal-lock.png
share/pixmaps/wicd/good-signal.png
share/pixmaps/wicd/high-signal-lock.png
share/pixmaps/wicd/high-signal.png
share/pixmaps/wicd/idle-bad-signal-lock.png
share/pixmaps/wicd/idle-bad-signal.png
share/pixmaps/wicd/idle-good-signal-lock.png
share/pixmaps/wicd/idle-good-signal.png
share/pixmaps/wicd/idle-high-signal-lock.png
share/pixmaps/wicd/idle-high-signal.png
share/pixmaps/wicd/idle-low-signal-lock.png
share/pixmaps/wicd/idle-low-signal.png
share/pixmaps/wicd/low-signal-lock.png
share/pixmaps/wicd/low-signal.png
share/pixmaps/wicd/no-signal.png
share/pixmaps/wicd/receiving-bad-signal-lock.png
share/pixmaps/wicd/receiving-bad-signal.png
share/pixmaps/wicd/receiving-good-signal-lock.png
share/pixmaps/wicd/receiving-good-signal.png
share/pixmaps/wicd/receiving-high-signal-lock.png
share/pixmaps/wicd/receiving-high-signal.png
share/pixmaps/wicd/receiving-low-signal-lock.png
share/pixmaps/wicd/receiving-low-signal.png
share/pixmaps/wicd/signal-100.png
share/pixmaps/wicd/signal-25.png
share/pixmaps/wicd/signal-50.png
share/pixmaps/wicd/signal-75.png
share/pixmaps/wicd/transmitting-bad-signal-lock.png
share/pixmaps/wicd/transmitting-bad-signal.png
share/pixmaps/wicd/transmitting-good-signal-lock.png
share/pixmaps/wicd/transmitting-good-signal.png
share/pixmaps/wicd/transmitting-high-signal-lock.png
share/pixmaps/wicd/transmitting-high-signal.png
share/pixmaps/wicd/transmitting-low-signal-lock.png
share/pixmaps/wicd/transmitting-low-signal.png
share/pixmaps/wicd/wired-gui.svg
share/pixmaps/wicd/wired.png
share/examples/rc.d/wicd
share/applications/wicd.desktop
@unexec ${RM} %D/etc/xdg/autostart/wicd-tray.desktop
@unexec ${RM} %D/etc/dbus-1/system.d/wicd.conf
@unexec ${RM} /var/lib/wicd/WHEREAREMYFILES
@unexec ${RMDIR} /var/lib/wicd/configurations && ${RMDIR} /var/lib/wicd
@pkgdir etc/wicd/scripts/predisconnect
@pkgdir etc/wicd/scripts/preconnect
@pkgdir etc/wicd/scripts/postdisconnect
@pkgdir etc/wicd/scripts/postconnect
