# py-jabber-yahoo-transport/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/10/05 03:34:01 schnoebe Exp $
${TRANSPORTDIR}/adhoc.py
${TRANSPORTDIR}/adhoc.pyc
${TRANSPORTDIR}/avatar.py
${TRANSPORTDIR}/avatar.pyc
${TRANSPORTDIR}/config.py
${TRANSPORTDIR}/config.pyc
${TRANSPORTDIR}/curphoo/YahooMD5.py
${TRANSPORTDIR}/curphoo/YahooMD5.pyc
${TRANSPORTDIR}/curphoo/__init__.py
${TRANSPORTDIR}/curphoo/__init__.pyc
${TRANSPORTDIR}/curphoo/cpformat.py
${TRANSPORTDIR}/curphoo/cpformat.pyc
${TRANSPORTDIR}/curphoo/md5crypt.py
${TRANSPORTDIR}/curphoo/md5crypt.pyc
${TRANSPORTDIR}/curphoo/pysha.py
${TRANSPORTDIR}/curphoo/pysha.pyc
${TRANSPORTDIR}/curphoo/yahoo_fn.py
${TRANSPORTDIR}/curphoo/yahoo_fn.pyc
${TRANSPORTDIR}/htmlutils.py
${TRANSPORTDIR}/htmlutils.pyc
${TRANSPORTDIR}/jep0133.py
${TRANSPORTDIR}/jep0133.pyc
${TRANSPORTDIR}/migratespool.py
${TRANSPORTDIR}/migratespool.pyc
${TRANSPORTDIR}/roomlist.py
${TRANSPORTDIR}/roomlist.pyc
${TRANSPORTDIR}/toolbox.py
${TRANSPORTDIR}/toolbox.pyc
${TRANSPORTDIR}/xmlconfig.py
${TRANSPORTDIR}/xmlconfig.pyc
${TRANSPORTDIR}/yahoo.py
${TRANSPORTDIR}/yahoo.pyc
${TRANSPORTDIR}/yahoo_helpers.py
${TRANSPORTDIR}/yahoo_helpers.pyc
${TRANSPORTDIR}/ylib.py
${TRANSPORTDIR}/ylib.pyc
${JABBER_EXDIR}/yahoo-transport.xml
share/examples/rc.d/yahoo-transport
