# smtp-gated/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/09/25 08:53:14 bartoszkuzma Exp $
man/man5/smtp-gated.conf.5
man/man8/smtp-gated.8
sbin/smtp-gated
share/doc/smtp-gated/README.PL
share/examples/rc.d/smtp_gated
share/examples/smtp-gated/smtp-gated.conf
@pkgdir share/doc/smtp-gated
@pkgdir share/examples/smtp-gated
