# pdnsd/PLIST
@comment $NetBSD: PLIST,v 1.2 2007/03/03 22:07:35 cheusov Exp $
man/man5/pdnsd.conf.5
man/man8/pdnsd-ctl.8
man/man8/pdnsd.8
sbin/pdnsd
sbin/pdnsd-ctl
share/examples/pdnsd/pdnsd.conf.sample
share/examples/rc.d/pdnsd
@dirrm share/examples/pdnsd
