# xinetd/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:34:07 jeremy-c-reed Exp $
man/man5/xinetd.conf.5
man/man8/itox.8
man/man8/xconv.pl.8
man/man8/xinetd.8
man/man8/xinetd.log.8
sbin/itox
sbin/xconv.pl
sbin/xinetd
share/examples/rc.d/xinetd
share/examples/xinetd/xinetd.conf
share/examples/xinetd/xinetd.d/chargen-dgram
share/examples/xinetd/xinetd.d/chargen-stream
share/examples/xinetd/xinetd.d/daytime-dgram
share/examples/xinetd/xinetd.d/daytime-stream
share/examples/xinetd/xinetd.d/discard-dgram
share/examples/xinetd/xinetd.d/discard-stream
share/examples/xinetd/xinetd.d/echo-dgram
share/examples/xinetd/xinetd.d/echo-stream
share/examples/xinetd/xinetd.d/ftp-sensor
share/examples/xinetd/xinetd.d/tcpmux-server
share/examples/xinetd/xinetd.d/time-dgram
share/examples/xinetd/xinetd.d/time-stream
@dirrm share/examples/xinetd/xinetd.d
@dirrm share/examples/xinetd
