# feathercoin/PLIST
@comment $NetBSD$
sbin/feathercoind
share/examples/feathercoin/feathercoin.conf
share/examples/rc.d/feathercoind
