# cvsnt/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/10/11 10:44:32 thomasklausner Exp $
bin/cvs
bin/cvslockd
bin/cvsnt
${AUDIT_MYSQL}
lib/cvsnt/mdns/mini.la
lib/cvsnt/protocols/enum.la
lib/cvsnt/protocols/ext.la
lib/cvsnt/protocols/pserver.la
lib/cvsnt/protocols/server.la
lib/cvsnt/protocols/sserver.la
lib/cvsnt/protocols/sspi.la
lib/cvsnt/triggers/audit.la
lib/cvsnt/triggers/email.la
lib/cvsnt/triggers/checkout.la
lib/cvsnt/triggers/info.la
lib/cvsnt/xdiff/extdiff.la
lib/libcvsapi.la
lib/libcvstools.la
lib/libmdnsclient.la
man/man1/cvs.1
man/man5/cvs.5
share/examples/rc.d/cvslockd
share/examples/cvsnt/ca.pem.example
share/examples/cvsnt/Plugins.example
share/examples/cvsnt/PServer.example
