# rsyslog/PLIST
@comment $NetBSD: PLIST,v 1.2 2008/08/10 23:22:05 mstapelberg Exp $
man/man5/rsyslog.conf.5
man/man8/rsyslogd.8
sbin/rsyslogd
share/doc/rsyslog/bugs.html
share/doc/rsyslog/features.html
share/doc/rsyslog/history.html
share/doc/rsyslog/how2help.html
share/doc/rsyslog/install.html
share/doc/rsyslog/manual.html
share/doc/rsyslog/property_replacer.html
share/doc/rsyslog/rsyslog_conf.html
share/doc/rsyslog/rsyslog_mysql.html
share/doc/rsyslog/rsyslog_php_syslog_ng.html
share/doc/rsyslog/rsyslog_stunnel.html
share/doc/rsyslog/status.html
share/doc/rsyslog/version_naming.html
share/examples/rc.d/rsyslogd
share/examples/rsyslog/createDB.sql
share/examples/rsyslog/rsyslog.conf
@dirrm share/doc/rsyslog
@dirrm share/examples/rsyslog
lib/rsyslog/imklog.la
lib/rsyslog/imklog.so
lib/rsyslog/immark.la
lib/rsyslog/immark.so
lib/rsyslog/imtcp.la
lib/rsyslog/imtcp.so
lib/rsyslog/imudp.la
lib/rsyslog/imudp.so
lib/rsyslog/imuxsock.la
lib/rsyslog/imuxsock.so
lib/rsyslog/lmnet.la
lib/rsyslog/lmnet.so
lib/rsyslog/lmnetstrms.la
lib/rsyslog/lmnetstrms.so
lib/rsyslog/lmnsd_gtls.la
lib/rsyslog/lmnsd_gtls.so
lib/rsyslog/lmnsd_ptcp.la
lib/rsyslog/lmnsd_ptcp.so
lib/rsyslog/lmregexp.la
lib/rsyslog/lmregexp.so
lib/rsyslog/lmtcpclt.la
lib/rsyslog/lmtcpclt.so
lib/rsyslog/lmtcpsrv.la
lib/rsyslog/lmtcpsrv.so
lib/rsyslog/omtesting.la
lib/rsyslog/omtesting.so
