# fail2ban/PLIST
@comment $NetBSD: PLIST,v 1.3 2012/04/08 20:22:11 thomasklausner Exp $
bin/fail2ban-client
bin/fail2ban-regex
bin/fail2ban-server
${PYSITELIB}/client/__init__.py
${PYSITELIB}/client/__init__.pyc
${PYSITELIB}/client/__init__.pyo
${PYSITELIB}/client/actionreader.py
${PYSITELIB}/client/actionreader.pyc
${PYSITELIB}/client/actionreader.pyo
${PYSITELIB}/client/beautifier.py
${PYSITELIB}/client/beautifier.pyc
${PYSITELIB}/client/beautifier.pyo
${PYSITELIB}/client/configparserinc.py
${PYSITELIB}/client/configparserinc.pyc
${PYSITELIB}/client/configparserinc.pyo
${PYSITELIB}/client/configreader.py
${PYSITELIB}/client/configreader.pyc
${PYSITELIB}/client/configreader.pyo
${PYSITELIB}/client/configurator.py
${PYSITELIB}/client/configurator.pyc
${PYSITELIB}/client/configurator.pyo
${PYSITELIB}/client/csocket.py
${PYSITELIB}/client/csocket.pyc
${PYSITELIB}/client/csocket.pyo
${PYSITELIB}/client/fail2banreader.py
${PYSITELIB}/client/fail2banreader.pyc
${PYSITELIB}/client/fail2banreader.pyo
${PYSITELIB}/client/filterreader.py
${PYSITELIB}/client/filterreader.pyc
${PYSITELIB}/client/filterreader.pyo
${PYSITELIB}/client/jailreader.py
${PYSITELIB}/client/jailreader.pyc
${PYSITELIB}/client/jailreader.pyo
${PYSITELIB}/client/jailsreader.py
${PYSITELIB}/client/jailsreader.pyc
${PYSITELIB}/client/jailsreader.pyo
${PYSITELIB}/common/__init__.py
${PYSITELIB}/common/__init__.pyc
${PYSITELIB}/common/__init__.pyo
${PYSITELIB}/common/helpers.py
${PYSITELIB}/common/helpers.pyc
${PYSITELIB}/common/helpers.pyo
${PYSITELIB}/common/protocol.py
${PYSITELIB}/common/protocol.pyc
${PYSITELIB}/common/protocol.pyo
${PYSITELIB}/common/version.py
${PYSITELIB}/common/version.pyc
${PYSITELIB}/common/version.pyo
${PYSITELIB}/${EGG_FILE}
${PYSITELIB}/server/__init__.py
${PYSITELIB}/server/__init__.pyc
${PYSITELIB}/server/__init__.pyo
${PYSITELIB}/server/action.py
${PYSITELIB}/server/action.pyc
${PYSITELIB}/server/action.pyo
${PYSITELIB}/server/actions.py
${PYSITELIB}/server/actions.pyc
${PYSITELIB}/server/actions.pyo
${PYSITELIB}/server/asyncserver.py
${PYSITELIB}/server/asyncserver.pyc
${PYSITELIB}/server/asyncserver.pyo
${PYSITELIB}/server/banmanager.py
${PYSITELIB}/server/banmanager.pyc
${PYSITELIB}/server/banmanager.pyo
${PYSITELIB}/server/datedetector.py
${PYSITELIB}/server/datedetector.pyc
${PYSITELIB}/server/datedetector.pyo
${PYSITELIB}/server/datetemplate.py
${PYSITELIB}/server/datetemplate.pyc
${PYSITELIB}/server/datetemplate.pyo
${PYSITELIB}/server/faildata.py
${PYSITELIB}/server/faildata.pyc
${PYSITELIB}/server/faildata.pyo
${PYSITELIB}/server/failmanager.py
${PYSITELIB}/server/failmanager.pyc
${PYSITELIB}/server/failmanager.pyo
${PYSITELIB}/server/failregex.py
${PYSITELIB}/server/failregex.pyc
${PYSITELIB}/server/failregex.pyo
${PYSITELIB}/server/filter.py
${PYSITELIB}/server/filter.pyc
${PYSITELIB}/server/filter.pyo
${PYSITELIB}/server/filtergamin.py
${PYSITELIB}/server/filtergamin.pyc
${PYSITELIB}/server/filtergamin.pyo
${PYSITELIB}/server/filterpoll.py
${PYSITELIB}/server/filterpoll.pyc
${PYSITELIB}/server/filterpoll.pyo
${PYSITELIB}/server/iso8601.py
${PYSITELIB}/server/iso8601.pyc
${PYSITELIB}/server/iso8601.pyo
${PYSITELIB}/server/jail.py
${PYSITELIB}/server/jail.pyc
${PYSITELIB}/server/jail.pyo
${PYSITELIB}/server/jails.py
${PYSITELIB}/server/jails.pyc
${PYSITELIB}/server/jails.pyo
${PYSITELIB}/server/jailthread.py
${PYSITELIB}/server/jailthread.pyc
${PYSITELIB}/server/jailthread.pyo
${PYSITELIB}/server/mytime.py
${PYSITELIB}/server/mytime.pyc
${PYSITELIB}/server/mytime.pyo
${PYSITELIB}/server/server.py
${PYSITELIB}/server/server.pyc
${PYSITELIB}/server/server.pyo
${PYSITELIB}/server/ticket.py
${PYSITELIB}/server/ticket.pyc
${PYSITELIB}/server/ticket.pyo
${PYSITELIB}/server/transmitter.py
${PYSITELIB}/server/transmitter.pyc
${PYSITELIB}/server/transmitter.pyo
man/man1/fail2ban-client.1
man/man1/fail2ban-regex.1
man/man1/fail2ban-server.1
share/examples/rc.d/fail2ban
