# hiawatha/PLIST
@comment $NetBSD$
bin/ssi-cgi
include/mbedtls/aes.h
include/mbedtls/aria.h
include/mbedtls/asn1.h
include/mbedtls/asn1write.h
include/mbedtls/base64.h
include/mbedtls/bignum.h
include/mbedtls/build_info.h
include/mbedtls/camellia.h
include/mbedtls/ccm.h
include/mbedtls/chacha20.h
include/mbedtls/chachapoly.h
include/mbedtls/check_config.h
include/mbedtls/cipher.h
include/mbedtls/cmac.h
include/mbedtls/compat-2.x.h
include/mbedtls/config_adjust_legacy_crypto.h
include/mbedtls/config_adjust_legacy_from_psa.h
include/mbedtls/config_adjust_psa_from_legacy.h
include/mbedtls/config_adjust_psa_superset_legacy.h
include/mbedtls/config_adjust_ssl.h
include/mbedtls/config_adjust_x509.h
include/mbedtls/config_psa.h
include/mbedtls/constant_time.h
include/mbedtls/ctr_drbg.h
include/mbedtls/debug.h
include/mbedtls/des.h
include/mbedtls/dhm.h
include/mbedtls/ecdh.h
include/mbedtls/ecdsa.h
include/mbedtls/ecjpake.h
include/mbedtls/ecp.h
include/mbedtls/entropy.h
include/mbedtls/error.h
include/mbedtls/gcm.h
include/mbedtls/hkdf.h
include/mbedtls/hmac_drbg.h
include/mbedtls/lms.h
include/mbedtls/mbedtls_config.h
include/mbedtls/md.h
include/mbedtls/md5.h
include/mbedtls/memory_buffer_alloc.h
include/mbedtls/net_sockets.h
include/mbedtls/nist_kw.h
include/mbedtls/oid.h
include/mbedtls/pem.h
include/mbedtls/pk.h
include/mbedtls/pkcs12.h
include/mbedtls/pkcs5.h
include/mbedtls/pkcs7.h
include/mbedtls/platform.h
include/mbedtls/platform_time.h
include/mbedtls/platform_util.h
include/mbedtls/poly1305.h
include/mbedtls/private_access.h
include/mbedtls/psa_util.h
include/mbedtls/ripemd160.h
include/mbedtls/rsa.h
include/mbedtls/sha1.h
include/mbedtls/sha256.h
include/mbedtls/sha3.h
include/mbedtls/sha512.h
include/mbedtls/ssl.h
include/mbedtls/ssl_cache.h
include/mbedtls/ssl_ciphersuites.h
include/mbedtls/ssl_cookie.h
include/mbedtls/ssl_ticket.h
include/mbedtls/threading.h
include/mbedtls/timing.h
include/mbedtls/version.h
include/mbedtls/x509.h
include/mbedtls/x509_crl.h
include/mbedtls/x509_crt.h
include/mbedtls/x509_csr.h
include/psa/build_info.h
include/psa/crypto.h
include/psa/crypto_adjust_auto_enabled.h
include/psa/crypto_adjust_config_key_pair_types.h
include/psa/crypto_adjust_config_synonyms.h
include/psa/crypto_builtin_composites.h
include/psa/crypto_builtin_key_derivation.h
include/psa/crypto_builtin_primitives.h
include/psa/crypto_compat.h
include/psa/crypto_config.h
include/psa/crypto_driver_common.h
include/psa/crypto_driver_contexts_composites.h
include/psa/crypto_driver_contexts_key_derivation.h
include/psa/crypto_driver_contexts_primitives.h
include/psa/crypto_extra.h
include/psa/crypto_legacy.h
include/psa/crypto_platform.h
include/psa/crypto_se_driver.h
include/psa/crypto_sizes.h
include/psa/crypto_struct.h
include/psa/crypto_types.h
include/psa/crypto_values.h
lib/hiawatha/letsencrypt/acmev2.php
lib/hiawatha/letsencrypt/config.php
lib/hiawatha/letsencrypt/hiawatha_config.php
lib/hiawatha/letsencrypt/http.php
lib/hiawatha/letsencrypt/https.php
lib/hiawatha/letsencrypt/letsencrypt.conf
lib/hiawatha/letsencrypt/letsencrypt.php
lib/hiawatha/letsencrypt/logfile.php
lib/hiawatha/letsencrypt/openssl.conf
lib/hiawatha/letsencrypt/rsa.php
lib/hiawatha/libmbedcrypto.a
lib/hiawatha/libmbedtls.a
lib/hiawatha/libmbedx509.a
man/man1/cgi-wrapper.1
man/man1/hiawatha.1
man/man1/lefh.1
man/man1/ssi-cgi.1
man/man1/wigwam.1
sbin/cgi-wrapper
sbin/hiawatha
sbin/lefh
sbin/wigwam
share/examples/hiawatha/cgi-wrapper.conf
share/examples/hiawatha/error.xslt
share/examples/hiawatha/hiawatha.conf
share/examples/hiawatha/index.xslt
share/examples/hiawatha/mimetype.conf
share/examples/rc.d/hiawatha
var/hiawatha/index.html
@pkgdir var/run
@pkgdir var/log/hiawatha
@pkgdir var/db/hiawatha
@pkgdir etc/hiawatha
