# sabnzbdplus/files/sabnzbd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
  Assumes user=sabnzbd group=sabnzbd
  Assumes '/opt/local/bin/' is installation directory
  
  To enable/disable
    svcadm enable sabnzbd
    svcadm disable sabnzbd
  
  To check logs
    tail /var/svc/log/network-sabnzbd\:default.log
-->
 
<service_bundle type='manifest' name='sabnzbd'>
 
  <service
    name='network/sabnzbd'
    type='service'
    version='1'>
 
    <create_default_instance enabled='false' />
    <single_instance />
 
   <!-- Only start in muti-user mode -->
   <dependency name='multi-user'
      grouping='require_all'
      restart_on='none'
      type='service'>
      <service_fmri value='svc:/milestone/multi-user' />
   </dependency>
 
    <!-- Wait for network interfaces to be initialized. -->
    <dependency name='network'
        grouping='require_all'
        restart_on='none'
        type='service'>
        <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
 
    <!-- Wait for all local filesystems to be mounted. -->
    <dependency name='filesystem-local'
        grouping='require_all'
        restart_on='none'
        type='service'>
        <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
 
    <!-- execute as user sabnzbd -->
    <method_context>
        <method_credential user='sabnzbd' group='sabnzbd' />
    </method_context>
 
    <exec_method
        type='method'
        name='start'
		exec='/opt/local/bin/SABnzbd.py -d -f /opt/local/sabnzbd/sabnzbd.ini'
        timeout_seconds='60'>
        <method_context>
	        <method_environment>
		        <envvar name='PATH' value='/opt/local/bin:/opt/local/sbin:/usr/bin:/usr/sbin
		        '/>
	        </method_environment>
        </method_context>
    </exec_method>
 
    <exec_method
        type='method'
        name='stop'
        exec=':kill'
        timeout_seconds='2'>
    </exec_method>
 
    <template>
       <common_name>
          <loctext xml:lang='C'>sabnzbd</loctext>
       </common_name>
       <documentation>
          <doc_link name='sabnzbd'
                    uri='www.http://sabnzbd.org/' />
       </documentation>
    </template>
 
  </service>
 
</service_bundle>

