# primecoin/PLIST
@comment $NetBSD$
sbin/primecoind
share/examples/primecoin/primecoin.conf
share/examples/rc.d/primecoind
