# pts/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/10/11 10:45:06 thomasklausner Exp $
sbin/pts
share/examples/rc.d/pts
share/pts/locale/en.msg
share/pts/locale/ja.msg
share/pts/pts.ban
share/pts/pts.ini
share/pts/pts.motd
share/pts/pts.news
share/pts/pts.secure
@exec ${MKDIR} %D/share/pts/logs
@exec ${MKDIR} %D/share/pts/dstats
