# pts/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:03 jeremy-c-reed Exp $
sbin/pts
share/examples/rc.d/pts
share/pts/locale/en.msg
share/pts/locale/ja.msg
share/pts/pts.ban
share/pts/pts.ini
share/pts/pts.motd
share/pts/pts.news
share/pts/pts.secure
@exec ${MKDIR} %D/share/pts/logs
@dirrm share/pts/logs
@dirrm share/pts/locale
@exec ${MKDIR} %D/share/pts/dstats
@dirrm share/pts/dstats
@dirrm share/pts
