# parpd/PLIST
@comment $NetBSD: PLIST,v 1.1 2008/10/07 21:12:14 rsmarples Exp $
man/man5/parpd.conf.5
man/man8/parpd.8
sbin/parpd
share/examples/rc.d/parpd
