# entrance/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/10/11 10:44:34 thomasklausner Exp $
bin/entrance
bin/entrance_edit
bin/entrance_edit-config
bin/entrance_wrapper
include/Entrance_Edit.h
lib/libentrance_edit.la
libexec/entrace/entrance_login
sbin/entranced
share/entrace/images/pointer.png
share/entrace/images/sessions/blackbox.png
share/entrace/images/sessions/default.png
share/entrace/images/sessions/enlightenment.png
share/entrace/images/sessions/failsafe.png
share/entrace/images/sessions/gnome.png
share/entrace/images/sessions/kde.png
share/entrace/images/sessions/xfce.png
share/entrace/themes/darkrock.edj
share/entrace/themes/default.edj
share/entrace/users/default.edj
share/examples/rc.d/entranced
