# am-utils/PLIST
@comment $NetBSD$
bin/pawd
info/am-utils.info
info/am-utils.info-1
info/am-utils.info-2
lib/libamu.la
man/man1/pawd.1
man/man5/amd.conf.5
man/man8/amd.8
man/man8/amq.8
man/man8/automount2amd.8
man/man8/fixmount.8
man/man8/fsinfo.8
man/man8/hlfsd.8
man/man8/mk-amd-map.8
man/man8/sun2amd.8
man/man8/wire-test.8
sbin/am-eject
sbin/amd
sbin/amd2ldif
sbin/amd2sun
sbin/amq
sbin/automount2amd
sbin/ctl-amd
sbin/ctl-hlfsd
sbin/fix-amd-map
sbin/fixmount
sbin/fixrmtab
sbin/fsinfo
sbin/hlfsd
sbin/lostaltmail
sbin/mk-amd-map
sbin/sun2amd
sbin/wait4amd
sbin/wait4amd2die
sbin/wire-test
share/doc/am-utils/AUTHORS
share/doc/am-utils/BUGS
share/doc/am-utils/COPYING
share/doc/am-utils/ChangeLog
share/doc/am-utils/FAQ
share/doc/am-utils/INSTALL
share/doc/am-utils/NEWS
share/doc/am-utils/README
share/doc/am-utils/README.attrcache
share/doc/am-utils/README.autofs
share/doc/am-utils/README.ldap
share/doc/am-utils/README.osx
share/doc/am-utils/README.y2k
share/examples/am-utils/amd.conf
share/examples/am-utils/lostaltmail.conf
share/examples/rc.d/amd-pkgsrc
