# ezbounce/PLIST
@comment $NetBSD: PLIST,v 1.6 2013/01/20 22:12:46 othyro Exp $
man/man1/ezbounce.1
sbin/ezbchk
sbin/ezbounce
share/doc/ezbounce/CHANGES
share/doc/ezbounce/README
share/doc/ezbounce/TODO
share/doc/ezbounce/motd.txt
share/doc/ezbounce/sample.conf
share/examples/ezbounce/ezbounce.conf.default
share/examples/rc.d/ezbounce
