# ezbounce/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:33:56 jeremy-c-reed Exp $
man/man1/ezbounce.1
sbin/ezbchk
sbin/ezbounce
sbin/mdidentd
share/doc/ezbounce/CHANGES
share/doc/ezbounce/README
share/doc/ezbounce/README.SSL
share/doc/ezbounce/TODO
share/doc/ezbounce/sample.MOTD
share/doc/ezbounce/sample.conf
share/examples/ezbounce/ezbounce.conf.default
share/examples/rc.d/ezbounce
@dirrm share/examples/ezbounce
@dirrm share/doc/ezbounce
