# kapacitor/files/smf/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='@SMF_NAME@:default'>
  <service name='@SMF_PREFIX@/@SMF_NAME@' type='service' version='1'>
    <create_default_instance enabled='false' />
    <single_instance />
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default' />
    </dependency>
    <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default' />
    </dependency>
    <method_context working_directory="@DATADIR@">
      <method_credential user='@KAPACITOR_USER@' group='@KAPACITOR_GROUP@' />
    </method_context>
    <exec_method type='method' name='start'
      exec='@PREFIX@/bin/kapacitord -config @PKG_SYSCONFDIR@/kapacitor.conf -pidfile @DATADIR@/kapacitord.pid 2>> @LOGDIR@/kapacitor.log &amp;'
      timeout_seconds="60" />
    <exec_method type='method' name='stop' exec=':kill' timeout_seconds="60" />
    <template>
      <common_name>
        <loctext xml:lang='C'>Kapacitor metrics monitoring agent</loctext>
      </common_name>
      <documentation>
        <manpage title='kapacitord' section='1' manpath='/@PREFIX@/man' />
        <doc_link name='homepage' uri='https://docs.influxdata.com/kapacitor/' />
      </documentation>
    </template>
  </service>
</service_bundle>
