# ejabberd/PLIST
@comment $NetBSD: PLIST,v 1.5 2008/02/14 12:21:07 sebpierrel Exp $
sbin/ejabberdctl
share/examples/rc.d/ejabberd
${EJABBERD_EXDIR}/ejabberd.cfg
${EJABBERD_EXDIR}/ejabberd.defaults
${EJABBERD_EXDIR}/ejabberd.inetrc
var/lib/ejabberd/ebin/ELDAPv3.beam
var/lib/ejabberd/ebin/XmppAddr.beam
var/lib/ejabberd/ebin/acl.beam
var/lib/ejabberd/ebin/adhoc.beam
var/lib/ejabberd/ebin/cyrsasl.beam
var/lib/ejabberd/ebin/cyrsasl_anonymous.beam
var/lib/ejabberd/ebin/cyrsasl_digest.beam
var/lib/ejabberd/ebin/cyrsasl_plain.beam
var/lib/ejabberd/ebin/ejabberd.app
var/lib/ejabberd/ebin/ejabberd.beam
var/lib/ejabberd/ebin/ejabberd_admin.beam
var/lib/ejabberd/ebin/ejabberd_app.beam
var/lib/ejabberd/ebin/ejabberd_auth.beam
var/lib/ejabberd/ebin/ejabberd_auth_anonymous.beam
var/lib/ejabberd/ebin/ejabberd_auth_external.beam
var/lib/ejabberd/ebin/ejabberd_auth_internal.beam
var/lib/ejabberd/ebin/ejabberd_auth_ldap.beam
var/lib/ejabberd/ebin/ejabberd_auth_odbc.beam
var/lib/ejabberd/ebin/ejabberd_auth_pam.beam
var/lib/ejabberd/ebin/ejabberd_c2s.beam
var/lib/ejabberd/ebin/ejabberd_c2s_config.beam
var/lib/ejabberd/ebin/ejabberd_config.beam
var/lib/ejabberd/ebin/ejabberd_ctl.beam
var/lib/ejabberd/ebin/ejabberd_frontend_socket.beam
var/lib/ejabberd/ebin/ejabberd_hooks.beam
var/lib/ejabberd/ebin/ejabberd_http.beam
var/lib/ejabberd/ebin/ejabberd_http_bind.beam
var/lib/ejabberd/ebin/ejabberd_http_poll.beam
var/lib/ejabberd/ebin/ejabberd_listener.beam
var/lib/ejabberd/ebin/ejabberd_local.beam
var/lib/ejabberd/ebin/ejabberd_logger_h.beam
var/lib/ejabberd/ebin/ejabberd_loglevel.beam
var/lib/ejabberd/ebin/ejabberd_node_groups.beam
var/lib/ejabberd/ebin/ejabberd_odbc.beam
var/lib/ejabberd/ebin/ejabberd_odbc_sup.beam
var/lib/ejabberd/ebin/ejabberd_rdbms.beam
var/lib/ejabberd/ebin/ejabberd_receiver.beam
var/lib/ejabberd/ebin/ejabberd_router.beam
var/lib/ejabberd/ebin/ejabberd_s2s.beam
var/lib/ejabberd/ebin/ejabberd_s2s_in.beam
var/lib/ejabberd/ebin/ejabberd_s2s_out.beam
var/lib/ejabberd/ebin/ejabberd_service.beam
var/lib/ejabberd/ebin/ejabberd_sm.beam
var/lib/ejabberd/ebin/ejabberd_socket.beam
var/lib/ejabberd/ebin/ejabberd_sup.beam
var/lib/ejabberd/ebin/ejabberd_system_monitor.beam
var/lib/ejabberd/ebin/ejabberd_tmp_sup.beam
var/lib/ejabberd/ebin/ejabberd_update.beam
var/lib/ejabberd/ebin/ejabberd_web.beam
var/lib/ejabberd/ebin/ejabberd_web_admin.beam
var/lib/ejabberd/ebin/ejabberd_zlib.beam
var/lib/ejabberd/ebin/ejd2odbc.beam
var/lib/ejabberd/ebin/eldap.beam
var/lib/ejabberd/ebin/eldap_filter.beam
var/lib/ejabberd/ebin/eldap_pool.beam
var/lib/ejabberd/ebin/eldap_utils.beam
var/lib/ejabberd/ebin/extauth.beam
var/lib/ejabberd/ebin/gen_iq_handler.beam
var/lib/ejabberd/ebin/gen_mod.beam
var/lib/ejabberd/ebin/gen_pubsub_node.beam
var/lib/ejabberd/ebin/gen_pubsub_nodetree.beam
var/lib/ejabberd/ebin/iconv.beam
var/lib/ejabberd/ebin/idna.beam
var/lib/ejabberd/ebin/jd2ejd.beam
var/lib/ejabberd/ebin/jlib.beam
var/lib/ejabberd/ebin/mod_adhoc.beam
var/lib/ejabberd/ebin/mod_announce.beam
var/lib/ejabberd/ebin/mod_caps.beam
var/lib/ejabberd/ebin/mod_configure.beam
var/lib/ejabberd/ebin/mod_configure2.beam
var/lib/ejabberd/ebin/mod_disco.beam
var/lib/ejabberd/ebin/mod_echo.beam
var/lib/ejabberd/ebin/mod_http_bind.beam
var/lib/ejabberd/ebin/mod_http_fileserver.beam
var/lib/ejabberd/ebin/mod_irc.beam
var/lib/ejabberd/ebin/mod_irc_connection.beam
var/lib/ejabberd/ebin/mod_last.beam
var/lib/ejabberd/ebin/mod_last_odbc.beam
var/lib/ejabberd/ebin/mod_muc.beam
var/lib/ejabberd/ebin/mod_muc_log.beam
var/lib/ejabberd/ebin/mod_muc_room.beam
var/lib/ejabberd/ebin/mod_offline.beam
var/lib/ejabberd/ebin/mod_offline_odbc.beam
var/lib/ejabberd/ebin/mod_privacy.beam
var/lib/ejabberd/ebin/mod_privacy_odbc.beam
var/lib/ejabberd/ebin/mod_private.beam
var/lib/ejabberd/ebin/mod_private_odbc.beam
var/lib/ejabberd/ebin/mod_proxy65.beam
var/lib/ejabberd/ebin/mod_proxy65_lib.beam
var/lib/ejabberd/ebin/mod_proxy65_service.beam
var/lib/ejabberd/ebin/mod_proxy65_sm.beam
var/lib/ejabberd/ebin/mod_proxy65_stream.beam
var/lib/ejabberd/ebin/mod_pubsub.beam
var/lib/ejabberd/ebin/mod_register.beam
var/lib/ejabberd/ebin/mod_roster.beam
var/lib/ejabberd/ebin/mod_roster_odbc.beam
var/lib/ejabberd/ebin/mod_service_log.beam
var/lib/ejabberd/ebin/mod_shared_roster.beam
var/lib/ejabberd/ebin/mod_stats.beam
var/lib/ejabberd/ebin/mod_time.beam
var/lib/ejabberd/ebin/mod_vcard.beam
var/lib/ejabberd/ebin/mod_vcard_ldap.beam
var/lib/ejabberd/ebin/mod_vcard_odbc.beam
var/lib/ejabberd/ebin/mod_version.beam
var/lib/ejabberd/ebin/node_buddy.beam
var/lib/ejabberd/ebin/node_club.beam
var/lib/ejabberd/ebin/node_default.beam
var/lib/ejabberd/ebin/node_dispatch.beam
var/lib/ejabberd/ebin/node_pep.beam
var/lib/ejabberd/ebin/node_private.beam
var/lib/ejabberd/ebin/node_public.beam
var/lib/ejabberd/ebin/nodetree_default.beam
var/lib/ejabberd/ebin/nodetree_virtual.beam
var/lib/ejabberd/ebin/odbc_queries.beam
var/lib/ejabberd/ebin/p1_fsm.beam
var/lib/ejabberd/ebin/p1_mnesia.beam
var/lib/ejabberd/ebin/ram_file_io_server.beam
var/lib/ejabberd/ebin/randoms.beam
var/lib/ejabberd/ebin/sha.beam
var/lib/ejabberd/ebin/shaper.beam
var/lib/ejabberd/ebin/stringprep.beam
var/lib/ejabberd/ebin/stringprep_sup.beam
var/lib/ejabberd/ebin/tls.beam
var/lib/ejabberd/ebin/translate.beam
var/lib/ejabberd/ebin/xml.beam
var/lib/ejabberd/ebin/xml_stream.beam
var/lib/ejabberd/priv/lib/ejabberd_zlib_drv.so
var/lib/ejabberd/priv/lib/expat_erl.so
var/lib/ejabberd/priv/lib/iconv_erl.so
var/lib/ejabberd/priv/lib/stringprep_drv.so
var/lib/ejabberd/priv/lib/tls_drv.so
var/lib/ejabberd/priv/msgs/ca.msg
var/lib/ejabberd/priv/msgs/cs.msg
var/lib/ejabberd/priv/msgs/de.msg
var/lib/ejabberd/priv/msgs/es.msg
var/lib/ejabberd/priv/msgs/fr.msg
var/lib/ejabberd/priv/msgs/gl.msg
var/lib/ejabberd/priv/msgs/it.msg
var/lib/ejabberd/priv/msgs/ja.msg
var/lib/ejabberd/priv/msgs/nl.msg
var/lib/ejabberd/priv/msgs/pl.msg
var/lib/ejabberd/priv/msgs/pt-br.msg
var/lib/ejabberd/priv/msgs/pt.msg
var/lib/ejabberd/priv/msgs/ru.msg
var/lib/ejabberd/priv/msgs/sk.msg
var/lib/ejabberd/priv/msgs/sv.msg
var/lib/ejabberd/priv/msgs/th.msg
var/lib/ejabberd/priv/msgs/tr.msg
var/lib/ejabberd/priv/msgs/uk.msg
var/lib/ejabberd/priv/msgs/vi.msg
var/lib/ejabberd/priv/msgs/wa.msg
var/lib/ejabberd/priv/msgs/zh.msg
share/doc/${DISTNAME}/Makefile
share/doc/${DISTNAME}/COPYING
share/doc/${DISTNAME}/dev.html
share/doc/${DISTNAME}/dev.tex
share/doc/${DISTNAME}/disco.png
share/doc/${DISTNAME}/discorus.png
share/doc/${DISTNAME}/features.html
share/doc/${DISTNAME}/features.tex
share/doc/${DISTNAME}/flow.dot
share/doc/${DISTNAME}/guide.html
share/doc/${DISTNAME}/guide.tex
share/doc/${DISTNAME}/guide.pdf
share/doc/${DISTNAME}/introduction.tex
share/doc/${DISTNAME}/logo.png
share/doc/${DISTNAME}/mod_http_bind.tex
share/doc/${DISTNAME}/mod_http_fileserver.tex
share/doc/${DISTNAME}/release_notes_0.9.1.txt
share/doc/${DISTNAME}/release_notes_0.9.8.txt
share/doc/${DISTNAME}/release_notes_0.9.txt
share/doc/${DISTNAME}/release_notes_1.0.0.txt
share/doc/${DISTNAME}/release_notes_1.1.0.txt
share/doc/${DISTNAME}/release_notes_1.1.1.txt
share/doc/${DISTNAME}/release_notes_1.1.2.txt
share/doc/${DISTNAME}/release_notes_1.1.3.txt
share/doc/${DISTNAME}/release_notes_1.1.4.txt
share/doc/${DISTNAME}/release_notes_2.0.0.txt
share/doc/${DISTNAME}/version.tex
share/doc/${DISTNAME}/webadmmain.png
share/doc/${DISTNAME}/webadmmainru.png
share/doc/${DISTNAME}/yozhikheader.png
share/doc/${DISTNAME}/api/Makefile
share/doc/${DISTNAME}/api/process-one.css
share/doc/${DISTNAME}/api/overview.edoc
@dirrm share/doc/${DISTNAME}/api
@dirrm share/doc/${DISTNAME}
@dirrm var/lib/ejabberd/priv/msgs
@dirrm var/lib/ejabberd/priv/lib
@dirrm var/lib/ejabberd/priv
@dirrm var/lib/ejabberd/ebin
@dirrm var/lib/ejabberd
