# shaper/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/10/11 10:45:17 thomasklausner Exp $
sbin/shapecfg
share/doc/shaper/README.shaper
share/examples/rc.d/shaper
share/examples/shaper/shaper.conf
