# xrdp09/PLIST
@comment $NetBSD$
bin/xrdp-dis
bin/xrdp-genkeymap
bin/xrdp-keygen
bin/xrdp-sesadmin
bin/xrdp-sesrun
include/ms-erref.h
include/ms-fscc.h
include/ms-rdpbcgr.h
include/ms-rdpeclip.h
include/ms-rdpedisp.h
include/ms-rdpefs.h
include/ms-rdpegdi.h
include/ms-rdpele.h
include/ms-rdperp.h
include/ms-smb2.h
include/xrdp_client_info.h
include/xrdp_constants.h
include/xrdp_rail.h
include/xrdp_sockets.h
lib/pkgconfig/xrdp.pc
lib/xrdp/libcommon.la
lib/xrdp/libmc.la
lib/xrdp/libscp.la
lib/xrdp/libvnc.la
lib/xrdp/libxrdp.la
lib/xrdp/libxrdpapi.la
lib/xrdp/libxup.la
man/man1/xrdp-dis.1
man/man5/sesman.ini.5
man/man5/xrdp.ini.5
man/man8/xrdp-chansrv.8
man/man8/xrdp-genkeymap.8
man/man8/xrdp-keygen.8
man/man8/xrdp-sesadmin.8
man/man8/xrdp-sesman.8
man/man8/xrdp-sesrun.8
man/man8/xrdp.8
sbin/xrdp
sbin/xrdp-chansrv
sbin/xrdp-sesman
share/examples/rc.d/sesman
share/examples/rc.d/xrdp
share/xrdp/ad24b.bmp
share/xrdp/ad256.bmp
share/xrdp/cursor0.cur
share/xrdp/cursor1.cur
share/xrdp/sans-10.fv1
share/xrdp/xrdp24b.bmp
share/xrdp/xrdp256.bmp
share/xrdp/xrdp_logo.bmp
