# up-imapproxy/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:07 jeremy-c-reed Exp $
sbin/in.imapproxyd
sbin/pimpstat
share/doc/up-imapproxy/README
share/doc/up-imapproxy/README.known_issues
share/doc/up-imapproxy/README.ssl
share/examples/rc.d/upimapproxy
share/examples/up-imapproxy/imapproxy.conf
@dirrm share/examples/up-imapproxy
@dirrm share/doc/up-imapproxy
