# ldapdns/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/10/11 10:44:51 thomasklausner Exp $
sbin/${PKGNAME}-dg
sbin/${PKGNAME}-hosts
sbin/${PKGNAME}-vc
share/doc/${PKGNAME}/INSTALL
share/doc/${PKGNAME}/README.axfr
share/doc/${PKGNAME}/README.charset
share/doc/${PKGNAME}/README.generic-rr
share/examples/rc.d/ldapdns
