# aolserver/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:33:38 jeremy-c-reed Exp $
share/examples/rc.d/nsd
