# pts/files/pts
#!/bin/sh

# PROVIDE: pts
# REQUIRE: DAEMON LOGIN network

. /etc/rc.subr

name="pts"
rcvar=$name
procname="perl:"
command="/usr/pkg/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/usr/pkg/share/pts/pts.ini"

load_rc_config $name
run_rc_command "$1"
