# riak/files/smf/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="riak">
  <service name="@SMF_PREFIX@/riak" type="service" version="1">
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <dependency name="epmd" grouping="require_any" restart_on="restart" type="service">
      <service_fmri value="svc:/@SMF_PREFIX@/riak-epmd:default" />
    </dependency>
    <method_context working_directory="@RIAK_DATA@" project="riak">
      <method_credential user="@RIAK_USER@" group="@RIAK_GROUP@" />
      <method_environment>
        <envvar name="HOME" value="@RIAK_DATA@" />
        <envvar name="LOGNAME" value="@RIAK_USER@" />
        <envvar name="PATH" value="/usr/bin:/bin:@PREFIX@/bin:@PREFIX@/sbin" />
        <envvar name="LD_PRELOAD_32" value="/lib/libumem.so.1" />
        <envvar name="LD_PRELOAD_64" value="/lib/64/libumem.so.1" />
        <envvar name="UMEM_OPTIONS" value="allocator=best" />
      </method_environment>
    </method_context>
    <exec_method type="method" name="start" exec="@PREFIX@/sbin/riak start" timeout_seconds="60" />
    <exec_method type="method" name="stop" exec="@PREFIX@/sbin/riak stop" timeout_seconds="60" />
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="contract" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>
    <stability value="Evolving" />
    <template>
      <common_name>
        <loctext xml:lang="C">Riak Data Store</loctext>
      </common_name>
    </template>
  </service>
  <service name="@SMF_PREFIX@/riak-epmd" type="service" version="1">
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <dependency name="erlang-epmd" grouping="exclude_all" restart_on="restart" type="service">
      <service_fmri value="svc:/@SMF_PREFIX@/epmd:default" />
    </dependency>
    <method_context>
      <method_credential user="@RIAK_USER@" group="@RIAK_GROUP@" />
    </method_context>
    <exec_method type="method" name="start" exec="@PREFIX@/riak/erts-@ERTS_VERSION@/bin/epmd -daemon" timeout_seconds="60" />
    <exec_method type="method" name="stop" exec="@PREFIX@/riak/erts-@ERTS_VERSION@/bin/epmd -kill" timeout_seconds="60" />
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="contract" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>
    <property_group name="application" type="application"></property_group>
    <stability value="Evolving" />
    <template>
      <common_name>
        <loctext xml:lang="C">Erlang Port Mapper Daemon</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
