# py-jabber-icqt/PLIST
@comment $NetBSD: PLIST,v 1.5 2006/12/20 01:06:57 schnoebe Exp $
bin/PyICQt
${TRANSPORTDIR}/src/avatar.py
${TRANSPORTDIR}/src/adhoc.py
${TRANSPORTDIR}/src/config.py
${TRANSPORTDIR}/src/bool.py
${TRANSPORTDIR}/src/langs/__init__.py
${TRANSPORTDIR}/src/langs/cs.py
${TRANSPORTDIR}/src/langs/de.py
${TRANSPORTDIR}/src/langs/en.py
${TRANSPORTDIR}/src/langs/es.py
${TRANSPORTDIR}/src/langs/fr.py
${TRANSPORTDIR}/src/langs/nl.py
${TRANSPORTDIR}/src/langs/pl.py
${TRANSPORTDIR}/src/langs/ru.py
${TRANSPORTDIR}/src/langs/sv.py
${TRANSPORTDIR}/src/utils.py
${TRANSPORTDIR}/src/contact.py
${TRANSPORTDIR}/src/debug.py
${TRANSPORTDIR}/src/disco.py
${TRANSPORTDIR}/src/globals.py
${TRANSPORTDIR}/src/imgmanip.py
${TRANSPORTDIR}/src/iq.py
${TRANSPORTDIR}/src/jabw.py
${TRANSPORTDIR}/src/lang.py
${TRANSPORTDIR}/src/legacy/services/ConfirmAccount.py
${TRANSPORTDIR}/src/legacy/services/EmailLookup.py
${TRANSPORTDIR}/src/legacy/services/__init__.py
${TRANSPORTDIR}/src/legacy/__init__.py
${TRANSPORTDIR}/src/legacy/buddies.py
${TRANSPORTDIR}/src/legacy/glue.py
${TRANSPORTDIR}/src/legacy/icqt.py
${TRANSPORTDIR}/src/pubsub.py
${TRANSPORTDIR}/src/main.py
${TRANSPORTDIR}/src/services/GatewayTranslator.py
${TRANSPORTDIR}/src/services/ConnectUsers.py
${TRANSPORTDIR}/src/services/EntityTime.py
${TRANSPORTDIR}/src/services/IqAvatarFactory.py
${TRANSPORTDIR}/src/services/LastActivity.py
${TRANSPORTDIR}/src/services/PingService.py
${TRANSPORTDIR}/src/services/RosterRetrieval.py
${TRANSPORTDIR}/src/services/SearchFactory.py
${TRANSPORTDIR}/src/services/Statistics.py
${TRANSPORTDIR}/src/services/VCardFactory.py
${TRANSPORTDIR}/src/services/VersionTeller.py
${TRANSPORTDIR}/src/services/__init__.py
${TRANSPORTDIR}/src/register.py
${TRANSPORTDIR}/src/sasl.py
${TRANSPORTDIR}/src/tlib/countrycodes.py
${TRANSPORTDIR}/src/tlib/__init__.py
${TRANSPORTDIR}/src/tlib/httpcompat.py
${TRANSPORTDIR}/src/tlib/ifcompat.py
${TRANSPORTDIR}/src/tlib/langcodes.py
${TRANSPORTDIR}/src/tlib/oscar.py
${TRANSPORTDIR}/src/tlib/scheduler.py
${TRANSPORTDIR}/src/tlib/socks5.py
${TRANSPORTDIR}/src/tlib/sockserror.py
${TRANSPORTDIR}/src/web/__init__.py
${TRANSPORTDIR}/src/web/handler.py
${TRANSPORTDIR}/src/web/xmppcred.py
${TRANSPORTDIR}/src/session.py
${TRANSPORTDIR}/src/svninfo.py
${TRANSPORTDIR}/src/xdb/__init__.py
${TRANSPORTDIR}/src/xdb/mysql.py
${TRANSPORTDIR}/src/xdb/template.py
${TRANSPORTDIR}/src/xdb/legacyaimtransport.py
${TRANSPORTDIR}/src/xdb/legacyjittransport.py
${TRANSPORTDIR}/src/xdb/xmlfiles.py
${TRANSPORTDIR}/src/twistfix/__init__.py
${TRANSPORTDIR}/src/twistfix/README
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/xmpp_stringprep.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/__init__.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/client.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/component.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/jid.py
${TRANSPORTDIR}/src/twistfix/words/protocols/jabber/jstrports.py
${TRANSPORTDIR}/src/twistfix/words/protocols/__init__.py
${TRANSPORTDIR}/src/twistfix/words/__init__.py
${TRANSPORTDIR}/src/twistfix/words/xish/xpathparser.py
${TRANSPORTDIR}/src/twistfix/words/xish/__init__.py
${TRANSPORTDIR}/src/twistfix/words/xish/domish.py
${TRANSPORTDIR}/src/twistfix/words/xish/utility.py
${TRANSPORTDIR}/src/twistfix/words/xish/xmlstream.py
${TRANSPORTDIR}/src/twistfix/words/xish/xpath.py
${TRANSPORTDIR}/src/twistfix/version.py
${TRANSPORTDIR}/src/xmlconfig.py
${TRANSPORTDIR}/data/www/css/style.css
${TRANSPORTDIR}/data/www/images/pyicq.png
${TRANSPORTDIR}/data/www/template.html
${TRANSPORTDIR}/data/defaultAIMAvatar.png
${TRANSPORTDIR}/data/defaultICQAvatar.png
${TRANSPORTDIR}/tools/db-setup.mysql
${TRANSPORTDIR}/tools/Makefile
${TRANSPORTDIR}/tools/infodump.py
${TRANSPORTDIR}/tools/managessi.py
${TRANSPORTDIR}/tools/migrate.py
${JABBER_EXDIR}/pyicqt.xml
share/examples/rc.d/pyicqt
@dirrm ${TRANSPORTDIR}/tools
@dirrm ${TRANSPORTDIR}/src/xdb
@dirrm ${TRANSPORTDIR}/src/web
@dirrm ${TRANSPORTDIR}/src/twistfix/words/xish
@dirrm ${TRANSPORTDIR}/src/twistfix/words/protocols/jabber
@dirrm ${TRANSPORTDIR}/src/twistfix/words/protocols
@dirrm ${TRANSPORTDIR}/src/twistfix/words
@dirrm ${TRANSPORTDIR}/src/twistfix
@dirrm ${TRANSPORTDIR}/src/tlib
@dirrm ${TRANSPORTDIR}/src/services
@dirrm ${TRANSPORTDIR}/src/legacy/services
@dirrm ${TRANSPORTDIR}/src/legacy
@dirrm ${TRANSPORTDIR}/src/langs
@dirrm ${TRANSPORTDIR}/src
@dirrm ${TRANSPORTDIR}/data/www/images
@dirrm ${TRANSPORTDIR}/data/www/css
@dirrm ${TRANSPORTDIR}/data/www
@dirrm ${TRANSPORTDIR}/data
@dirrm ${TRANSPORTDIR}
