# py-jabber-icqt/PLIST
@comment $NetBSD: PLIST,v 1.9 2009/10/11 10:45:10 thomasklausner Exp $
bin/PyICQt
${JABBER_EXDIR}/pyicqt.xml
${TRANSPORTDIR}/data/defaultAIMAvatar.png
${TRANSPORTDIR}/data/defaultICQAvatar.png
${TRANSPORTDIR}/data/www/css/style.css
${TRANSPORTDIR}/data/www/images/pyicq.png
${TRANSPORTDIR}/data/www/template.html
${TRANSPORTDIR}/src/adhoc.py
${TRANSPORTDIR}/src/avatar.py
${TRANSPORTDIR}/src/bool.py
${TRANSPORTDIR}/src/chardet_utf/__init__.py
${TRANSPORTDIR}/src/chardet_utf/charsetgroupprober.py
${TRANSPORTDIR}/src/chardet_utf/charsetprober.py
${TRANSPORTDIR}/src/chardet_utf/codingstatemachine.py
${TRANSPORTDIR}/src/chardet_utf/constants.py
${TRANSPORTDIR}/src/chardet_utf/mbcsgroupprober.py
${TRANSPORTDIR}/src/chardet_utf/mbcssm.py
${TRANSPORTDIR}/src/chardet_utf/universaldetector.py
${TRANSPORTDIR}/src/chardet_utf/utf16beprober.py
${TRANSPORTDIR}/src/chardet_utf/utf16leprober.py
${TRANSPORTDIR}/src/chardet_utf/utf8prober.py
${TRANSPORTDIR}/src/config.py
${TRANSPORTDIR}/src/contact.py
${TRANSPORTDIR}/src/debug.py
${TRANSPORTDIR}/src/disco.py
${TRANSPORTDIR}/src/git_ver.py
${TRANSPORTDIR}/src/globals.py
${TRANSPORTDIR}/src/imgmanip.py
${TRANSPORTDIR}/src/iq.py
${TRANSPORTDIR}/src/jabw.py
${TRANSPORTDIR}/src/lang.py
${TRANSPORTDIR}/src/langs/__init__.py
${TRANSPORTDIR}/src/langs/cs.py
${TRANSPORTDIR}/src/langs/de.py
${TRANSPORTDIR}/src/langs/en.py
${TRANSPORTDIR}/src/langs/es.py
${TRANSPORTDIR}/src/langs/fr.py
${TRANSPORTDIR}/src/langs/nl.py
${TRANSPORTDIR}/src/langs/pl.py
${TRANSPORTDIR}/src/langs/ru.py
${TRANSPORTDIR}/src/langs/sv.py
${TRANSPORTDIR}/src/legacy/__init__.py
${TRANSPORTDIR}/src/legacy/buddies.py
${TRANSPORTDIR}/src/legacy/glue.py
${TRANSPORTDIR}/src/legacy/icqt.py
${TRANSPORTDIR}/src/legacy/services/ConfirmAccount.py
${TRANSPORTDIR}/src/legacy/services/EmailLookup.py
${TRANSPORTDIR}/src/legacy/services/Help.py
${TRANSPORTDIR}/src/legacy/services/SetXStatus.py
${TRANSPORTDIR}/src/legacy/services/Settings.py
${TRANSPORTDIR}/src/legacy/services/__init__.py
${TRANSPORTDIR}/src/main.py
${TRANSPORTDIR}/src/pubsub.py
${TRANSPORTDIR}/src/register.py
${TRANSPORTDIR}/src/sasl.py
${TRANSPORTDIR}/src/services/ConnectUsers.py
${TRANSPORTDIR}/src/services/EntityTime.py
${TRANSPORTDIR}/src/services/GatewayTranslator.py
${TRANSPORTDIR}/src/services/IqAvatarFactory.py
${TRANSPORTDIR}/src/services/LastActivity.py
${TRANSPORTDIR}/src/services/PingService.py
${TRANSPORTDIR}/src/services/RosterRetrieval.py
${TRANSPORTDIR}/src/services/SearchFactory.py
${TRANSPORTDIR}/src/services/Statistics.py
${TRANSPORTDIR}/src/services/VCardFactory.py
${TRANSPORTDIR}/src/services/VersionTeller.py
${TRANSPORTDIR}/src/services/XMPPPing.py
${TRANSPORTDIR}/src/services/__init__.py
${TRANSPORTDIR}/src/session.py
${TRANSPORTDIR}/src/git_ver.py
${TRANSPORTDIR}/src/tlib/__init__.py
${TRANSPORTDIR}/src/tlib/countrycodes.py
${TRANSPORTDIR}/src/tlib/httpcompat.py
${TRANSPORTDIR}/src/tlib/ifcompat.py
${TRANSPORTDIR}/src/tlib/langcodes.py
${TRANSPORTDIR}/src/tlib/oscar.py
${TRANSPORTDIR}/src/tlib/scheduler.py
${TRANSPORTDIR}/src/tlib/socks5.py
${TRANSPORTDIR}/src/tlib/sockserror.py
${TRANSPORTDIR}/src/twistd.py
${TRANSPORTDIR}/src/utils.py
${TRANSPORTDIR}/src/web/__init__.py
${TRANSPORTDIR}/src/web/handler.py
${TRANSPORTDIR}/src/web/xmppcred.py
${TRANSPORTDIR}/src/xdb/__init__.py
${TRANSPORTDIR}/src/xdb/legacyaimtransport.py
${TRANSPORTDIR}/src/xdb/legacyjittransport.py
${TRANSPORTDIR}/src/xdb/mysql.py
${TRANSPORTDIR}/src/xdb/template.py
${TRANSPORTDIR}/src/xdb/xmlfiles.py
${TRANSPORTDIR}/src/xmlconfig.py
${TRANSPORTDIR}/tools/Makefile
${TRANSPORTDIR}/tools/db-setup.mysql
${TRANSPORTDIR}/tools/infodump.py
${TRANSPORTDIR}/tools/managessi.py
${TRANSPORTDIR}/tools/memory_
${TRANSPORTDIR}/tools/migrate.py
${TRANSPORTDIR}/tools/threads_
share/examples/rc.d/pyicqt
