# hal-dev/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2007/07/22 10:46:04 krub1 Exp $
bin/hal-device
bin/hal-disable-polling
bin/hal-find-by-capability
bin/hal-find-by-property
bin/hal-get-property
bin/hal-is-caller-locked-out
bin/hal-lock
bin/hal-set-property
bin/lshal
include/hal/libhal-storage.h
include/hal/libhal.h
lib/hal/scripts/hal-functions
lib/hal/scripts/hal-luks-remove
lib/hal/scripts/hal-luks-setup
lib/hal/scripts/hal-luks-teardown
lib/hal/scripts/hal-system-killswitch-get-power
lib/hal/scripts/hal-system-killswitch-set-power
lib/hal/scripts/hal-system-lcd-get-brightness
lib/hal/scripts/hal-system-lcd-set-brightness
lib/hal/scripts/hal-system-power-hibernate
lib/hal/scripts/hal-system-power-reboot
lib/hal/scripts/hal-system-power-set-power-save
lib/hal/scripts/hal-system-power-shutdown
lib/hal/scripts/hal-system-power-suspend
lib/hal/scripts/hal-system-power-suspend-hybrid
lib/libhal-storage.la
lib/libhal.la
lib/pkgconfig/hal-storage.pc
lib/pkgconfig/hal.pc
libexec/hal-storage-cleanup-all-mountpoints
libexec/hal-storage-cleanup-mountpoint
libexec/hal-storage-closetray
libexec/hal-storage-eject
libexec/hal-storage-mount
libexec/hal-storage-unmount
libexec/hal-system-power-pm-is-supported
libexec/hal-system-power-pmu
libexec/hald-generate-fdi-cache
libexec/hald-runner
man/man1/hal-disable-polling.1
man/man1/hal-find-by-capability.1
man/man1/hal-find-by-property.1
man/man1/hal-get-property.1
man/man1/hal-is-caller-locked-out.1
man/man1/hal-lock.1
man/man1/hal-set-property.1
man/man1/lshal.1
man/man8/hald.8
sbin/hald
share/examples/rc.d/hald
share/hal/fdi/fdi.dtd
share/hal/fdi/policy/10osvendor/10-input-policy.fdi
share/hal/fdi/policy/10osvendor/10-laptop-panel-mgmt-policy.fdi
share/hal/fdi/policy/10osvendor/10-power-mgmt-policy.fdi
share/hal/fdi/policy/10osvendor/10-rfkill-switch.fdi
share/hal/fdi/policy/10osvendor/15-storage-luks.fdi
share/hal/fdi/policy/10osvendor/20-storage-methods.fdi
@exec ${MKDIR} %D/share/hal/fdi/policy/20thirdparty
@dirrm share/hal/fdi/policy/20thirdparty
@dirrm share/hal/fdi/policy/10osvendor
@dirrm share/hal/fdi/policy
@dirrm share/hal/fdi
@dirrm lib/hal/scripts
@dirrm lib/hal
@dirrm include/hal
@exec ${MKDIR} %D/etc/hal/fdi/policy
@dirrm etc/hal/fdi/policy
@dirrm etc/hal/fdi
@dirrm etc/dbus-1/system.d
