# callweaver/PLIST
@comment $NetBSD: PLIST,v 1.6 2010/01/13 13:21:04 x400 Exp $
bin/callweaver
sbin/streamplayer
include/callweaver/acl.h
include/callweaver/adsi.h
include/callweaver/aes.h
include/callweaver/alaw.h
include/callweaver/app.h
include/callweaver/callweaver_db.h
include/callweaver/callweaver_expr.h
include/callweaver/callweaver_hash.h
include/callweaver/callweaver_keywords.h
include/callweaver/callweaver_mm.h
include/callweaver/causes.h
include/callweaver/cdr.h
include/callweaver/chan_capi.h
include/callweaver/chan_capi20.h
include/callweaver/channel.h
include/callweaver/chanvars.h
include/callweaver/cli.h
include/callweaver/compiler.h
include/callweaver/config.h
include/callweaver/crypto.h
include/callweaver/cwobj.h
include/callweaver/devicestate.h
include/callweaver/dlfcn-compat.h
include/callweaver/dns.h
include/callweaver/dnsmgr.h
include/callweaver/dsp.h
include/callweaver/dundi.h
include/callweaver/enum.h
include/callweaver/features.h
include/callweaver/file.h
include/callweaver/frame.h
include/callweaver/generator.h
include/callweaver/generic_jb.h
include/callweaver/image.h
include/callweaver/indications.h
include/callweaver/inline_api.h
include/callweaver/io.h
include/callweaver/linkedlists.h
include/callweaver/localtime.h
include/callweaver/lock.h
include/callweaver/logger.h
include/callweaver/manager.h
include/callweaver/module.h
include/callweaver/monitor.h
include/callweaver/musiconhold.h
include/callweaver/netsock.h
include/callweaver/ogi.h
include/callweaver/options.h
include/callweaver/pbx.h
include/callweaver/phone_no_utils.h
include/callweaver/poll-compat.h
include/callweaver/privacy.h
include/callweaver/res_odbc.h
include/callweaver/resonator.h
include/callweaver/rtp.h
include/callweaver/say.h
include/callweaver/sched.h
include/callweaver/slinfactory.h
include/callweaver/srv.h
include/callweaver/strings.h
include/callweaver/stun.h
include/callweaver/term.h
include/callweaver/time.h
include/callweaver/transcap.h
include/callweaver/translate.h
include/callweaver/udp.h
include/callweaver/udpfromto.h
include/callweaver/udptl.h
include/callweaver/ulaw.h
include/callweaver/unaligned.h
include/callweaver/utils.h
include/callweaver/xlaw.h
lib/callweaver/libcallweaver.la
lib/callweaver/libcwjb.la
lib/callweaver/libcwsqlite3.la
lib/callweaver/modules/app_adsiprog.la
lib/callweaver/modules/app_amd.la
lib/callweaver/modules/app_authenticate.la
lib/callweaver/modules/app_backgrounddetect.la
lib/callweaver/modules/app_backticks.la
lib/callweaver/modules/app_cdr.la
lib/callweaver/modules/app_changrab.la
lib/callweaver/modules/app_chanisavail.la
lib/callweaver/modules/app_chanspy.la
lib/callweaver/modules/app_controlplayback.la
lib/callweaver/modules/app_db.la
lib/callweaver/modules/app_devstate.la
lib/callweaver/modules/app_dial.la
lib/callweaver/modules/app_dictate.la
lib/callweaver/modules/app_directed_pickup.la
lib/callweaver/modules/app_directory.la
lib/callweaver/modules/app_disa.la
lib/callweaver/modules/app_dtmftotext.la
lib/callweaver/modules/app_dumpchan.la
lib/callweaver/modules/app_echo.la
lib/callweaver/modules/app_enumlookup.la
lib/callweaver/modules/app_eval.la
lib/callweaver/modules/app_exec.la
lib/callweaver/modules/app_faxdetect.la
lib/callweaver/modules/app_forkcdr.la
lib/callweaver/modules/app_getcpeid.la
lib/callweaver/modules/app_getdevstate.la
lib/callweaver/modules/app_getextstate.la
lib/callweaver/modules/app_groupcount.la
lib/callweaver/modules/app_hasnewvoicemail.la
lib/callweaver/modules/app_lookupblacklist.la
lib/callweaver/modules/app_lookupcidname.la
lib/callweaver/modules/app_milliwatt.la
lib/callweaver/modules/app_muxmon.la
lib/callweaver/modules/app_nconference.la
lib/callweaver/modules/app_page.la
lib/callweaver/modules/app_pipe.la
lib/callweaver/modules/app_playback.la
lib/callweaver/modules/app_privacy.la
lib/callweaver/modules/app_proc.la
lib/callweaver/modules/app_queue.la
lib/callweaver/modules/app_random.la
lib/callweaver/modules/app_read.la
lib/callweaver/modules/app_record.la
lib/callweaver/modules/app_rxfax.la
lib/callweaver/modules/app_sayunixtime.la
lib/callweaver/modules/app_senddtmf.la
lib/callweaver/modules/app_sendtext.la
lib/callweaver/modules/app_setcallerpres.la
lib/callweaver/modules/app_setcdruserfield.la
lib/callweaver/modules/app_setrdnis.la
lib/callweaver/modules/app_settransfercapability.la
lib/callweaver/modules/app_sms.la
lib/callweaver/modules/app_softhangup.la
lib/callweaver/modules/app_stack.la
lib/callweaver/modules/app_system.la
lib/callweaver/modules/app_t38gateway.la
lib/callweaver/modules/app_transfer.la
lib/callweaver/modules/app_txfax.la
lib/callweaver/modules/app_userevent.la
lib/callweaver/modules/app_valetparking.la
lib/callweaver/modules/app_verbose.la
lib/callweaver/modules/app_voicemail.la
lib/callweaver/modules/app_waitfordigits.la
lib/callweaver/modules/app_waitforring.la
lib/callweaver/modules/app_waitforsilence.la
lib/callweaver/modules/app_while.la
lib/callweaver/modules/app_zapateller.la
lib/callweaver/modules/cdr_csv.la
lib/callweaver/modules/cdr_custom.la
lib/callweaver/modules/cdr_manager.la
lib/callweaver/modules/cdr_sqlite3.la
lib/callweaver/modules/chan_agent.la
lib/callweaver/modules/chan_features.la
lib/callweaver/modules/chan_iax2.la
lib/callweaver/modules/chan_local.la
lib/callweaver/modules/chan_mgcp.la
lib/callweaver/modules/chan_sip.la
lib/callweaver/modules/chan_woomera.la
lib/callweaver/modules/codec_a_mu.la
lib/callweaver/modules/codec_alaw.la
lib/callweaver/modules/codec_dvi_adpcm.la
lib/callweaver/modules/codec_g722.la
lib/callweaver/modules/codec_g722_16k_8k.la
lib/callweaver/modules/codec_g726.la
lib/callweaver/modules/codec_gsm.la
lib/callweaver/modules/codec_lpc10.la
lib/callweaver/modules/codec_oki_adpcm.la
lib/callweaver/modules/codec_speex.la
lib/callweaver/modules/codec_ulaw.la
lib/callweaver/modules/format_au.la
lib/callweaver/modules/format_g723_1.la
lib/callweaver/modules/format_g726.la
lib/callweaver/modules/format_g729.la
lib/callweaver/modules/format_gsm.la
lib/callweaver/modules/format_h263.la
lib/callweaver/modules/format_jpeg.la
lib/callweaver/modules/format_pcm.la
lib/callweaver/modules/format_pcm_alaw.la
lib/callweaver/modules/format_sln.la
lib/callweaver/modules/format_wav.la
lib/callweaver/modules/format_wav_gsm.la
lib/callweaver/modules/func_callerid.la
lib/callweaver/modules/func_cdr.la
lib/callweaver/modules/func_config.la
lib/callweaver/modules/func_db.la
lib/callweaver/modules/func_enum.la
lib/callweaver/modules/func_env.la
lib/callweaver/modules/func_fileexists.la
lib/callweaver/modules/func_groupcount.la
lib/callweaver/modules/func_language.la
lib/callweaver/modules/func_logic.la
lib/callweaver/modules/func_math.la
lib/callweaver/modules/func_md5.la
lib/callweaver/modules/func_moh.la
lib/callweaver/modules/func_strings.la
lib/callweaver/modules/func_timeout.la
lib/callweaver/modules/func_uri.la
lib/callweaver/modules/pbx_config.la
lib/callweaver/modules/pbx_dundi.la
lib/callweaver/modules/pbx_loopback.la
lib/callweaver/modules/pbx_realtime.la
lib/callweaver/modules/pbx_spool.la
lib/callweaver/modules/res_adsi.la
lib/callweaver/modules/res_crypto.la
lib/callweaver/modules/res_features.la
lib/callweaver/modules/res_indications.la
lib/callweaver/modules/res_monitor.la
lib/callweaver/modules/res_musiconhold.la
lib/callweaver/modules/res_ogi.la
${PLIST.jabber}lib/callweaver/modules/res_jabber.la
${PLIST.mysql}lib/callweaver/modules/app_sql_mysql.la
${PLIST.mysql}lib/callweaver/modules/cdr_mysql.la
${PLIST.mysql}lib/callweaver/modules/res_config_mysql.la
${PLIST.pgsql}lib/callweaver/modules/app_sql_postgres.la
${PLIST.pgsql}lib/callweaver/modules/cdr_pgsql.la
${PLIST.pgsql}lib/callweaver/modules/res_config_pgsql.la
${PLIST.snmp}lib/callweaver/modules/res_snmp.la
man/man8/callweaver.8
man/man8/safe_callweaver.8
sbin/callweaver_cli
sbin/safe_callweaver
${PLIST.sipsecurity}share/callweaver/keys/dh1024.pem
${PLIST.sipsecurity}share/callweaver/keys/dh512.pem
${PLIST.sipsecurity}share/callweaver/keys/servercert.pem
${PLIST.sipsecurity}share/callweaver/keys/serverkey.pem
${PLIST.sipsecurity}share/callweaver/keys/trustcerts.pem
share/callweaver/ogi/eogi-sphinx-test
share/callweaver/ogi/eogi-test
share/callweaver/ogi/fastogi-test
share/callweaver/ogi/ogi-test.ogi
share/doc/callweaver/CREDITS
share/doc/callweaver/LICENSE
share/doc/callweaver/README
share/examples/callweaver/adsi.conf
share/examples/callweaver/adtranvofr.conf
share/examples/callweaver/agents.conf
share/examples/callweaver/amd.conf
share/examples/callweaver/callweaver.adsi
share/examples/callweaver/callweaver.conf
share/examples/callweaver/cdr.conf
share/examples/callweaver/cdr_custom.conf
share/examples/callweaver/cdr_manager.conf
share/examples/callweaver/cdr_tds.conf
share/examples/callweaver/codecs.conf
share/examples/callweaver/db-memcached.conf
share/examples/callweaver/dnsmgr.conf
share/examples/callweaver/dundi.conf
share/examples/callweaver/enum.conf
share/examples/callweaver/extconfig.conf
share/examples/callweaver/extensions.ael
share/examples/callweaver/extensions.conf
share/examples/callweaver/features.conf
share/examples/callweaver/iax.conf
share/examples/callweaver/indications.conf
share/examples/callweaver/logger.conf
share/examples/callweaver/manager.conf
share/examples/callweaver/meetme.conf
share/examples/callweaver/mgcp.conf
share/examples/callweaver/modem.conf
share/examples/callweaver/modules.conf
share/examples/callweaver/musiconhold.conf
share/examples/callweaver/muted.conf
share/examples/callweaver/osp.conf
share/examples/callweaver/privacy.conf
share/examples/callweaver/queues.conf
share/examples/callweaver/res_snmp.conf
share/examples/callweaver/rtp.conf
share/examples/callweaver/sip.conf
share/examples/callweaver/sip_notify.conf
share/examples/callweaver/udptl.conf
share/examples/callweaver/voicemail.conf
share/examples/callweaver/woomera.conf
${PLIST.jabber}share/examples/callweaver/res_jabber.conf
${PLIST.mysql}share/examples/callweaver/cdr_mysql.conf
${PLIST.mysql}share/examples/callweaver/res_mysql.conf
${PLIST.pgsql}share/examples/callweaver/cdr_pgsql.conf
${PLIST.pgsql}share/examples/callweaver/res_pgsql.conf
share/examples/rc.d/callweaver
@exec ${MKDIR} %D/share/callweaver/moh
@exec ${MKDIR} %D/share/callweaver/keys
@exec ${MKDIR} %D/share/callweaver/images
@exec ${MKDIR} %D/lib/callweaver/modules/icd
