# spamassassin-cvs/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:05 jeremy-c-reed Exp $
@comment Most of this PLIST is auto-generated; see Makefile.
share/doc/html/spamassassin/ninjabutton.png
@comment @dirrm share/doc/html/spamassassin
share/examples/rc.d/spamd
share/examples/spamassassin/local.cf
share/examples/spamassassin/netbsd_lists.cf
share/examples/spamassassin/procmailrc.example
share/spamassassin/languages
share/spamassassin/triplets.txt
share/spamassassin/user_prefs.template
@comment @dirrm share/examples/spamassassin
