# polipo/PLIST
@comment $NetBSD: PLIST,v 1.5 2012/02/18 21:31:43 ishit Exp $
bin/polipo
info/polipo.info
man/man1/polipo.1
share/examples/rc.d/polipo
share/polipo/README
share/polipo/config.sample
share/polipo/forbidden.sample
share/polipo/www/Access-control.html
share/polipo/www/Allowed-ports.html
share/polipo/www/Asynchronous-writing.html
share/polipo/www/Background.html
share/polipo/www/Browser-configuration.html
share/polipo/www/Cache-transparency.html
share/polipo/www/Caching.html
share/polipo/www/Censor-Accept_002dLanguage.html
share/polipo/www/Censoring-headers.html
share/polipo/www/Chunk-memory.html
share/polipo/www/Client-connections.html
share/polipo/www/Concept-index.html
share/polipo/www/Configuring-Polipo.html
share/polipo/www/Contacting-servers.html
share/polipo/www/Copying.html
share/polipo/www/DNS.html
share/polipo/www/Daemon.html
share/polipo/www/Disk-cache.html
share/polipo/www/Disk-format.html
share/polipo/www/External-redirectors.html
share/polipo/www/Forbidden-Tunnels.html
share/polipo/www/Forbidden.html
share/polipo/www/HTTP-parent-proxies.html
share/polipo/www/HTTP-tuning.html
share/polipo/www/Intermediate-proxies.html
share/polipo/www/Internal-forbidden-list.html
share/polipo/www/Latency-and-throughput.html
share/polipo/www/Limiting-chunk-usage.html
share/polipo/www/Limiting-memory-usage.html
share/polipo/www/Limiting-object-usage.html
share/polipo/www/Local-server.html
share/polipo/www/Logging.html
share/polipo/www/Malloc-memory.html
share/polipo/www/Memory-cache.html
share/polipo/www/Memory-usage.html
share/polipo/www/Modifying-the-on_002ddisk-cache.html
share/polipo/www/Network-traffic.html
share/polipo/www/Network.html
share/polipo/www/OS-usage-limits.html
share/polipo/www/Offline-browsing.html
share/polipo/www/PMM.html
share/polipo/www/POST-and-PUT.html
share/polipo/www/Parent-proxies.html
share/polipo/www/Partial-instances.html
share/polipo/www/Persistent-connections.html
share/polipo/www/Pipelining.html
share/polipo/www/Polipo-Invocation.html
share/polipo/www/Poor-Mans-Multiplexing.html
share/polipo/www/Proxies-and-caches.html
share/polipo/www/Purging.html
share/polipo/www/Running.html
share/polipo/www/SOCKS-parent-proxies.html
share/polipo/www/Server-statistics.html
share/polipo/www/Server_002dside-behaviour.html
share/polipo/www/Stopping.html
share/polipo/www/The-web.html
share/polipo/www/Tuning-POST-and-PUT.html
share/polipo/www/Tuning-the-HTTP-parser.html
share/polipo/www/Tuning-validation.html
share/polipo/www/Tunnelling-connections.html
share/polipo/www/Tweaking-validation.html
share/polipo/www/Variable-index.html
share/polipo/www/Web-interface.html
share/polipo/www/index.html
@exec ${MKDIR} %D/etc/polipo
