# pdns2/PLIST
@comment $NetBSD: PLIST,v 1.1 2006/05/04 20:31:53 anielski Exp $
bin/pdns_control
bin/zone2ldap
bin/zone2sql
man/man8/pdns_control.8
man/man8/pdns_server.8
man/man8/zone2sql.8
sbin/pdns_server
share/examples/rc.d/pdns
