# dancer-ircd/PLIST
@comment $NetBSD: PLIST,v 1.2 2003/05/10 08:48:39 xtraeme Exp $
etc/rc.d/dancer-ircd
man/man8/ircd.8
sbin/dancer-ircd
sbin/install_ircd
sbin/mkconf
sbin/mkpasswd
sbin/viconf
share/doc/dancer-ircd/CIDR.txt
share/doc/dancer-ircd/README
share/doc/dancer-ircd/README.FIRST.hybrid-6
share/doc/dancer-ircd/README.TSora
share/doc/dancer-ircd/README.small_nets
share/doc/dancer-ircd/README.umodes
share/doc/dancer-ircd/RELNOTES.hybrid-6
share/doc/dancer-ircd/blalloc.txt
share/doc/dancer-ircd/hybrid/ChangeLog.hybrid-6.gz
share/doc/dancer-ircd/hybrid/README
share/doc/dancer-ircd/hybrid/example.conf
share/doc/dancer-ircd/hybrid/example.conf.trillian
share/doc/dancer-ircd/hybrid/mtrie.txt
share/doc/dancer-ircd/hybrid/opermyth.txt
share/doc/dancer-ircd/hybrid/opers.txt
share/doc/dancer-ircd/hybrid/server-version-info
share/doc/dancer-ircd/hybrid/simple.conf
share/doc/dancer-ircd/rfc1459.txt
share/examples/dancer-ircd/ircd.conf.default
@dirrm share/examples/dancer-ircd
@dirrm share/doc/dancer-ircd/hybrid
@dirrm share/doc/dancer-ircd
