# tor/PLIST
@comment $NetBSD: PLIST,v 1.2 2007/07/30 22:47:36 netcap Exp $
bin/tor
bin/tor-gencert
bin/tor-resolve
bin/torify
man/man1/tor-resolve.1
man/man1/tor.1
man/man1/torify.1
share/examples/rc.d/tor
share/examples/tor/tor-tsocks.conf
share/examples/tor/torrc.sample
@dirrm share/examples/tor
@dirrm etc/tor
