# milter-greylist-devel/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2008/10/20 06:58:21 pettai Exp $
bin/milter-greylist
man/man8/milter-greylist.8
man/man5/greylist.conf.5
share/doc/milter-greylist/README
share/examples/milter-greylist/greylist.conf
share/examples/rc.d/milter-greylist
@dirrm share/examples/milter-greylist
@dirrm share/doc/milter-greylist
