# ndevd/PLIST
@comment $NetBSD$
include/ndevd.h
man/man8/ndevd.8
sbin/ndevd
share/examples/rc.d/ndevd
