# quarkcoin/PLIST
@comment $NetBSD: PLIST,v 1.1 2014/09/22 11:22:56 othyro Exp $
bin/quarkcoin-qt
sbin/quarkcoind
share/examples/quarkcoin/quarkcoin.conf
share/examples/rc.d/quarkcoind
