# py-albatross/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:03 jeremy-c-reed Exp $
bin/al-httpd
bin/al-session-daemon
${PYSITELIB}/albatross/__init__.py
${PYSITELIB}/albatross/__init__.pyc
${PYSITELIB}/albatross/__init__.pyo
${PYSITELIB}/albatross/apacheapp.py
${PYSITELIB}/albatross/apacheapp.pyc
${PYSITELIB}/albatross/apacheapp.pyo
${PYSITELIB}/albatross/app.py
${PYSITELIB}/albatross/app.pyc
${PYSITELIB}/albatross/app.pyo
${PYSITELIB}/albatross/cgiapp.py
${PYSITELIB}/albatross/cgiapp.pyc
${PYSITELIB}/albatross/cgiapp.pyo
${PYSITELIB}/albatross/common.py
${PYSITELIB}/albatross/common.pyc
${PYSITELIB}/albatross/common.pyo
${PYSITELIB}/albatross/context.py
${PYSITELIB}/albatross/context.pyc
${PYSITELIB}/albatross/context.pyo
${PYSITELIB}/albatross/fcgiapp.py
${PYSITELIB}/albatross/fcgiapp.pyc
${PYSITELIB}/albatross/fcgiapp.pyo
${PYSITELIB}/albatross/httpdapp.py
${PYSITELIB}/albatross/httpdapp.pyc
${PYSITELIB}/albatross/httpdapp.pyo
${PYSITELIB}/albatross/pidfile.py
${PYSITELIB}/albatross/pidfile.pyc
${PYSITELIB}/albatross/pidfile.pyo
${PYSITELIB}/albatross/randompage.py
${PYSITELIB}/albatross/randompage.pyc
${PYSITELIB}/albatross/randompage.pyo
${PYSITELIB}/albatross/session.py
${PYSITELIB}/albatross/session.pyc
${PYSITELIB}/albatross/session.pyo
${PYSITELIB}/albatross/sessionfile.py
${PYSITELIB}/albatross/sessionfile.pyc
${PYSITELIB}/albatross/sessionfile.pyo
${PYSITELIB}/albatross/simpleserver.py
${PYSITELIB}/albatross/simpleserver.pyc
${PYSITELIB}/albatross/simpleserver.pyo
${PYSITELIB}/albatross/tags.py
${PYSITELIB}/albatross/tags.pyc
${PYSITELIB}/albatross/tags.pyo
${PYSITELIB}/albatross/template.py
${PYSITELIB}/albatross/template.pyc
${PYSITELIB}/albatross/template.pyo
share/examples/albatross/extension/cal.html
share/examples/albatross/extension/cal.py
share/examples/albatross/extension/install.py
share/examples/albatross/form4/form.html
share/examples/albatross/form4/form.py
share/examples/albatross/form4/install.py
share/examples/albatross/images/close.png
share/examples/albatross/images/ellipsis.png
share/examples/albatross/images/install.py
share/examples/albatross/images/open.png
share/examples/albatross/install.py
share/examples/albatross/mpperf/htaccess
share/examples/albatross/mpperf/install.py
share/examples/albatross/mpperf/main.html
share/examples/albatross/mpperf/main.py
share/examples/albatross/mpperf/mpperf.py
share/examples/albatross/paginate/install.py
share/examples/albatross/paginate/paginate.html
share/examples/albatross/paginate/paginate.py
share/examples/albatross/popview1/detail.html
share/examples/albatross/popview1/install.py
share/examples/albatross/popview1/list.html
share/examples/albatross/popview1/login.html
share/examples/albatross/popview1/popview.py
share/examples/albatross/popview1/popviewlib.py
share/examples/albatross/popview2/detail.html
share/examples/albatross/popview2/install.py
share/examples/albatross/popview2/list.html
share/examples/albatross/popview2/login.html
share/examples/albatross/popview2/popview.py
share/examples/albatross/popview2/popviewlib.py
share/examples/albatross/popview3/detail.html
share/examples/albatross/popview3/install.py
share/examples/albatross/popview3/list.html
share/examples/albatross/popview3/login.html
share/examples/albatross/popview3/popview.py
share/examples/albatross/popview3/popviewlib.py
share/examples/albatross/popview4/detail.html
share/examples/albatross/popview4/detail.py
share/examples/albatross/popview4/install.py
share/examples/albatross/popview4/list.html
share/examples/albatross/popview4/list.py
share/examples/albatross/popview4/login.html
share/examples/albatross/popview4/login.py
share/examples/albatross/popview4/popview.py
share/examples/albatross/popview4/popviewlib.py
share/examples/albatross/popview5/detail.html
share/examples/albatross/popview5/detail.py
share/examples/albatross/popview5/htaccess
share/examples/albatross/popview5/install.py
share/examples/albatross/popview5/list.html
share/examples/albatross/popview5/list.py
share/examples/albatross/popview5/login.html
share/examples/albatross/popview5/login.py
share/examples/albatross/popview5/popview.py
share/examples/albatross/popview5/popviewlib.py
share/examples/albatross/random/install.py
share/examples/albatross/random/paginate.html
share/examples/albatross/random/paginate.py
share/examples/albatross/random/random.py
share/examples/albatross/random/tree.html
share/examples/albatross/random/tree.py
share/examples/albatross/sybase/htaccess
share/examples/albatross/sybase/install.py
share/examples/albatross/sybase/macros.html
share/examples/albatross/sybase/sybase.py
share/examples/albatross/sybase/table-desc.html
share/examples/albatross/sybase/table-list.html
share/examples/albatross/sybase/table-view.html
share/examples/albatross/sybase/traceback.html
share/examples/albatross/templates/content1/content.py
share/examples/albatross/templates/content1/install.py
share/examples/albatross/templates/content1/main.html
share/examples/albatross/templates/content1/oops.html
share/examples/albatross/templates/content1/other.html
share/examples/albatross/templates/content2/content.py
share/examples/albatross/templates/content2/install.py
share/examples/albatross/templates/content2/macros.html
share/examples/albatross/templates/content2/main.html
share/examples/albatross/templates/content2/oops.html
share/examples/albatross/templates/content2/other.html
share/examples/albatross/templates/form1/form-display.html
share/examples/albatross/templates/form1/form.html
share/examples/albatross/templates/form1/form.py
share/examples/albatross/templates/form1/install.py
share/examples/albatross/templates/form2/form-display.html
share/examples/albatross/templates/form2/form.html
share/examples/albatross/templates/form2/form.py
share/examples/albatross/templates/form2/install.py
share/examples/albatross/templates/form3/form-display.html
share/examples/albatross/templates/form3/form.html
share/examples/albatross/templates/form3/form.py
share/examples/albatross/templates/form3/install.py
share/examples/albatross/templates/simple1/install.py
share/examples/albatross/templates/simple1/simple.py
share/examples/albatross/templates/simple2/install.py
share/examples/albatross/templates/simple2/simple.html
share/examples/albatross/templates/simple2/simple.py
share/examples/albatross/templates/simple3/install.py
share/examples/albatross/templates/simple3/simple.html
share/examples/albatross/templates/simple3/simple.py
share/examples/albatross/templates/simple4/install.py
share/examples/albatross/templates/simple4/simple.html
share/examples/albatross/templates/simple4/simple.py
share/examples/albatross/templates/simple5/install.py
share/examples/albatross/templates/simple5/simple.html
share/examples/albatross/templates/simple5/simple.py
share/examples/albatross/templates/stream/install.py
share/examples/albatross/templates/stream/stream.html
share/examples/albatross/templates/stream/stream.py
share/examples/albatross/templates/tree/tree.html
share/examples/albatross/templates/tree/tree.py
share/examples/albatross/tree1/install.py
share/examples/albatross/tree1/tree.html
share/examples/albatross/tree1/tree.py
share/examples/albatross/tree2/install.py
share/examples/albatross/tree2/tree.html
share/examples/albatross/tree2/tree.py
share/examples/albatross/tree3/install.py
share/examples/albatross/tree3/tree.html
share/examples/albatross/tree3/tree.py
share/examples/rc.d/al_session_daemon
@dirrm share/examples/albatross/tree3
@dirrm share/examples/albatross/tree2
@dirrm share/examples/albatross/tree1
@dirrm share/examples/albatross/templates/tree
@dirrm share/examples/albatross/templates/stream
@dirrm share/examples/albatross/templates/simple5
@dirrm share/examples/albatross/templates/simple4
@dirrm share/examples/albatross/templates/simple3
@dirrm share/examples/albatross/templates/simple2
@dirrm share/examples/albatross/templates/simple1
@dirrm share/examples/albatross/templates/form3
@dirrm share/examples/albatross/templates/form2
@dirrm share/examples/albatross/templates/form1
@dirrm share/examples/albatross/templates/content2
@dirrm share/examples/albatross/templates/content1
@dirrm share/examples/albatross/templates
@dirrm share/examples/albatross/sybase
@dirrm share/examples/albatross/random
@dirrm share/examples/albatross/popview5
@dirrm share/examples/albatross/popview4
@dirrm share/examples/albatross/popview3
@dirrm share/examples/albatross/popview2
@dirrm share/examples/albatross/popview1
@dirrm share/examples/albatross/paginate
@dirrm share/examples/albatross/mpperf
@dirrm share/examples/albatross/images
@dirrm share/examples/albatross/form4
@dirrm share/examples/albatross/extension
@dirrm share/examples/albatross
@dirrm ${PYSITELIB}/albatross
