# ftpd-tls/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/10/11 10:44:37 thomasklausner Exp $
libexec/ftpd-tls
man/man8/ftpd-tls.8
share/doc/${PKGBASE}/README.TLS
share/examples/rc.d/ftpdtls
