# libvirt/PLIST.Linux
@comment $NetBSD$
bin/virt-login-shell
lib/libvirt/connection-driver/libvirt_driver_lxc.la
lib/libvirt/connection-driver/libvirt_driver_network.la
lib/libvirt/connection-driver/libvirt_driver_nwfilter.la
lib/libvirt/storage-backend/libvirt_storage_backend_logical.la
lib/libvirt/storage-backend/libvirt_storage_backend_mpath.la
lib/libvirt/storage-backend/libvirt_storage_backend_scsi.la
lib/sysctl.d/60-libvirtd.conf
lib/systemd/system/libvirt-guests.service
lib/systemd/system/libvirtd-admin.socket
lib/systemd/system/libvirtd-ro.socket
lib/systemd/system/libvirtd-tcp.socket
lib/systemd/system/libvirtd-tls.socket
lib/systemd/system/libvirtd.service
lib/systemd/system/libvirtd.socket
lib/systemd/system/virt-guest-shutdown.target
lib/systemd/system/virtlockd-admin.socket
lib/systemd/system/virtlockd.service
lib/systemd/system/virtlockd.socket
lib/systemd/system/virtlogd-admin.socket
lib/systemd/system/virtlogd.service
lib/systemd/system/virtlogd.socket
lib/systemd/system/virtlxcd-admin.socket
lib/systemd/system/virtlxcd-ro.socket
lib/systemd/system/virtlxcd.service
lib/systemd/system/virtlxcd.socket
lib/systemd/system/virtnetworkd-admin.socket
lib/systemd/system/virtnetworkd-ro.socket
lib/systemd/system/virtnetworkd.service
lib/systemd/system/virtnetworkd.socket
lib/systemd/system/virtnwfilterd-admin.socket
lib/systemd/system/virtnwfilterd-ro.socket
lib/systemd/system/virtnwfilterd.service
lib/systemd/system/virtnwfilterd.socket
lib/systemd/system/virtproxyd-admin.socket
lib/systemd/system/virtproxyd-ro.socket
lib/systemd/system/virtproxyd-tcp.socket
lib/systemd/system/virtproxyd-tls.socket
lib/systemd/system/virtproxyd.service
lib/systemd/system/virtproxyd.socket
lib/systemd/system/virtqemud-admin.socket
lib/systemd/system/virtqemud-ro.socket
lib/systemd/system/virtqemud.service
lib/systemd/system/virtqemud.socket
lib/systemd/system/virtsecretd-admin.socket
lib/systemd/system/virtsecretd-ro.socket
lib/systemd/system/virtsecretd.service
lib/systemd/system/virtsecretd.socket
lib/systemd/system/virtstoraged-admin.socket
lib/systemd/system/virtstoraged-ro.socket
lib/systemd/system/virtstoraged.service
lib/systemd/system/virtstoraged.socket
lib/systemd/system/virtvboxd-admin.socket
lib/systemd/system/virtvboxd-ro.socket
lib/systemd/system/virtvboxd.service
lib/systemd/system/virtvboxd.socket
libexec/libvirt_leaseshelper
libexec/libvirt_lxc
libexec/virt-login-shell-helper
man/man1/virt-login-shell.1
sbin/virtlxcd
sbin/virtnetworkd
sbin/virtnwfilterd
share/augeas/lenses/libvirtd_lxc.aug
share/augeas/lenses/tests/test_libvirtd_lxc.aug
share/augeas/lenses/tests/test_virtlxcd.aug
share/augeas/lenses/tests/test_virtnetworkd.aug
share/augeas/lenses/tests/test_virtnwfilterd.aug
share/augeas/lenses/virtlxcd.aug
share/augeas/lenses/virtnetworkd.aug
share/augeas/lenses/virtnwfilterd.aug
share/examples/libvirt/lxc.conf
share/examples/libvirt/nwfilter/allow-arp.xml
share/examples/libvirt/nwfilter/allow-dhcp-server.xml
share/examples/libvirt/nwfilter/allow-dhcp.xml
share/examples/libvirt/nwfilter/allow-incoming-ipv4.xml
share/examples/libvirt/nwfilter/allow-ipv4.xml
share/examples/libvirt/nwfilter/clean-traffic-gateway.xml
share/examples/libvirt/nwfilter/clean-traffic.xml
share/examples/libvirt/nwfilter/no-arp-ip-spoofing.xml
share/examples/libvirt/nwfilter/no-arp-mac-spoofing.xml
share/examples/libvirt/nwfilter/no-arp-spoofing.xml
share/examples/libvirt/nwfilter/no-ip-multicast.xml
share/examples/libvirt/nwfilter/no-ip-spoofing.xml
share/examples/libvirt/nwfilter/no-mac-broadcast.xml
share/examples/libvirt/nwfilter/no-mac-spoofing.xml
share/examples/libvirt/nwfilter/no-other-l2-traffic.xml
share/examples/libvirt/nwfilter/no-other-rarp-traffic.xml
share/examples/libvirt/nwfilter/qemu-announce-self-rarp.xml
share/examples/libvirt/nwfilter/qemu-announce-self.xml
share/examples/libvirt/qemu/networks/autostart/default.xml
share/examples/libvirt/qemu/networks/default.xml
share/examples/libvirt/virt-login-shell.conf
share/examples/libvirt/virtlxcd.conf
share/examples/libvirt/virtnetworkd.conf
share/examples/libvirt/virtnwfilterd.conf
share/examples/sysconfig/libvirt-guests
share/examples/sysconfig/libvirtd
share/examples/sysconfig/virtlockd
share/examples/sysconfig/virtlogd
