# mopher/PLIST
@comment $NetBSD: PLIST,v 1.2 2012/12/11 00:39:42 pebog Exp $
${PLIST.bdb}lib/mopher/bdb.so
${PLIST.mysql}lib/mopher/sakila.so
${PLIST.spf}lib/mopher/spf.so
lib/mopher/cast.so
lib/mopher/counter.so
lib/mopher/list.so
lib/mopher/memdb.so
lib/mopher/rbl.so
lib/mopher/regex.so
lib/mopher/spamd.so
lib/mopher/string.so
sbin/moco
sbin/mopherd
share/examples/mopher/mail.acl.example
share/examples/mopher/mopherd.conf.example
share/examples/rc.d/mopherd
