# syslog-ng/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:06 jeremy-c-reed Exp $
bin/syslog2ng
man/man5/syslog-ng.conf.5
man/man8/syslog-ng.8
sbin/syslog-ng
share/doc/html/syslog-ng/book1.html
share/doc/html/syslog-ng/destinations.html
share/doc/html/syslog-ng/intro.html
share/doc/html/syslog-ng/logpath.html
share/doc/html/syslog-ng/msgroute.html
share/doc/html/syslog-ng/reference.html
share/doc/html/syslog-ng/tuning.html
share/doc/html/syslog-ng/x172.html
share/doc/html/syslog-ng/x361.html
share/doc/html/syslog-ng/x728.html
share/doc/html/syslog-ng/x768.html
share/doc/html/syslog-ng/x900.html
share/doc/html/syslog-ng/x906.html
share/doc/html/syslog-ng/x97.html
share/doc/syslog-ng/INSTALL
share/doc/syslog-ng/README
share/doc/syslog-ng/syslog-ng.conf.doc
share/examples/rc.d/syslogng
share/examples/syslog-ng/syslog-ng.conf.demo
share/examples/syslog-ng/syslog-ng.conf.sample
share/examples/syslog-ng/syslog-ng.conf.solaris
@dirrm share/examples/syslog-ng
@dirrm share/doc/syslog-ng
@dirrm share/doc/html/syslog-ng
