# zetacoin/PLIST
@comment $NetBSD$
bin/zetacoin-cli
${PLIST.qt}bin/zetacoin-qt
sbin/zetacoind
share/examples/zetacoin/zetacoin.conf
share/examples/rc.d/zetacoind
