# mathopd/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:00 jeremy-c-reed Exp $
sbin/mathopd
share/doc/mathopd/CHANGES
share/doc/mathopd/COPYING
share/doc/mathopd/README
share/doc/mathopd/cgi.txt
share/doc/mathopd/config.txt
share/doc/mathopd/running.txt
share/doc/mathopd/sample.cfg
share/doc/mathopd/syntax.txt
share/examples/mathopd/mathopd.conf
share/examples/rc.d/mathopd
@dirrm share/examples/mathopd
@dirrm share/doc/mathopd
