# mathopd/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/10/11 10:44:56 thomasklausner Exp $
sbin/mathopd
share/doc/mathopd/CHANGES
share/doc/mathopd/COPYING
share/doc/mathopd/README
share/doc/mathopd/cgi.txt
share/doc/mathopd/config.txt
share/doc/mathopd/running.txt
share/doc/mathopd/sample.cfg
share/doc/mathopd/syntax.txt
share/examples/mathopd/mathopd.conf
share/examples/rc.d/mathopd
