# py-bravo/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2011/06/01 03:51:46 schnoebe Exp $
${PYSITELIB}/${EGG_INFODIR}/PKG-INFO
${PYSITELIB}/${EGG_INFODIR}/SOURCES.txt
${PYSITELIB}/${EGG_INFODIR}/dependency_links.txt
${PYSITELIB}/${EGG_INFODIR}/requires.txt
${PYSITELIB}/${EGG_INFODIR}/top_level.txt
${PYSITELIB}/bravo/__init__.py
${PYSITELIB}/bravo/__init__.pyc
${PYSITELIB}/bravo/__init__.pyo
${PYSITELIB}/bravo/amp.py
${PYSITELIB}/bravo/amp.pyc
${PYSITELIB}/bravo/amp.pyo
${PYSITELIB}/bravo/blocks.py
${PYSITELIB}/bravo/blocks.pyc
${PYSITELIB}/bravo/blocks.pyo
${PYSITELIB}/bravo/chunk.py
${PYSITELIB}/bravo/chunk.pyc
${PYSITELIB}/bravo/chunk.pyo
${PYSITELIB}/bravo/config.py
${PYSITELIB}/bravo/config.pyc
${PYSITELIB}/bravo/config.pyo
${PYSITELIB}/bravo/entity.py
${PYSITELIB}/bravo/entity.pyc
${PYSITELIB}/bravo/entity.pyo
${PYSITELIB}/bravo/errors.py
${PYSITELIB}/bravo/errors.pyc
${PYSITELIB}/bravo/errors.pyo
${PYSITELIB}/bravo/factories/__init__.py
${PYSITELIB}/bravo/factories/__init__.pyc
${PYSITELIB}/bravo/factories/__init__.pyo
${PYSITELIB}/bravo/factories/beta.py
${PYSITELIB}/bravo/factories/beta.pyc
${PYSITELIB}/bravo/factories/beta.pyo
${PYSITELIB}/bravo/factories/infini.py
${PYSITELIB}/bravo/factories/infini.pyc
${PYSITELIB}/bravo/factories/infini.pyo
${PYSITELIB}/bravo/ibravo.py
${PYSITELIB}/bravo/ibravo.pyc
${PYSITELIB}/bravo/ibravo.pyo
${PYSITELIB}/bravo/inventory.py
${PYSITELIB}/bravo/inventory.pyc
${PYSITELIB}/bravo/inventory.pyo
${PYSITELIB}/bravo/irc.py
${PYSITELIB}/bravo/irc.pyc
${PYSITELIB}/bravo/irc.pyo
${PYSITELIB}/bravo/location.py
${PYSITELIB}/bravo/location.pyc
${PYSITELIB}/bravo/location.pyo
${PYSITELIB}/bravo/motd.py
${PYSITELIB}/bravo/motd.pyc
${PYSITELIB}/bravo/motd.pyo
${PYSITELIB}/bravo/nbt.py
${PYSITELIB}/bravo/nbt.pyc
${PYSITELIB}/bravo/nbt.pyo
${PYSITELIB}/bravo/packets/__init__.py
${PYSITELIB}/bravo/packets/__init__.pyc
${PYSITELIB}/bravo/packets/__init__.pyo
${PYSITELIB}/bravo/packets/beta.py
${PYSITELIB}/bravo/packets/beta.pyc
${PYSITELIB}/bravo/packets/beta.pyo
${PYSITELIB}/bravo/packets/encodings.py
${PYSITELIB}/bravo/packets/encodings.pyc
${PYSITELIB}/bravo/packets/encodings.pyo
${PYSITELIB}/bravo/packets/infini.py
${PYSITELIB}/bravo/packets/infini.pyc
${PYSITELIB}/bravo/packets/infini.pyo
${PYSITELIB}/bravo/plugin.py
${PYSITELIB}/bravo/plugin.pyc
${PYSITELIB}/bravo/plugin.pyo
${PYSITELIB}/bravo/plugins/__init__.py
${PYSITELIB}/bravo/plugins/__init__.pyc
${PYSITELIB}/bravo/plugins/__init__.pyo
${PYSITELIB}/bravo/plugins/authenticators.py
${PYSITELIB}/bravo/plugins/authenticators.pyc
${PYSITELIB}/bravo/plugins/authenticators.pyo
${PYSITELIB}/bravo/plugins/automatons.py
${PYSITELIB}/bravo/plugins/automatons.pyc
${PYSITELIB}/bravo/plugins/automatons.pyo
${PYSITELIB}/bravo/plugins/build_hooks.py
${PYSITELIB}/bravo/plugins/build_hooks.pyc
${PYSITELIB}/bravo/plugins/build_hooks.pyo
${PYSITELIB}/bravo/plugins/commands.py
${PYSITELIB}/bravo/plugins/commands.pyc
${PYSITELIB}/bravo/plugins/commands.pyo
${PYSITELIB}/bravo/plugins/compound_hooks.py
${PYSITELIB}/bravo/plugins/compound_hooks.pyc
${PYSITELIB}/bravo/plugins/compound_hooks.pyo
${PYSITELIB}/bravo/plugins/debug.py
${PYSITELIB}/bravo/plugins/debug.pyc
${PYSITELIB}/bravo/plugins/debug.pyo
${PYSITELIB}/bravo/plugins/dig_hooks.py
${PYSITELIB}/bravo/plugins/dig_hooks.pyc
${PYSITELIB}/bravo/plugins/dig_hooks.pyo
${PYSITELIB}/bravo/plugins/generators.py
${PYSITELIB}/bravo/plugins/generators.pyc
${PYSITELIB}/bravo/plugins/generators.pyo
${PYSITELIB}/bravo/plugins/paintings.py
${PYSITELIB}/bravo/plugins/paintings.pyc
${PYSITELIB}/bravo/plugins/paintings.pyo
${PYSITELIB}/bravo/plugins/physics.py
${PYSITELIB}/bravo/plugins/physics.pyc
${PYSITELIB}/bravo/plugins/physics.pyo
${PYSITELIB}/bravo/plugins/recipes.py
${PYSITELIB}/bravo/plugins/recipes.pyc
${PYSITELIB}/bravo/plugins/recipes.pyo
${PYSITELIB}/bravo/plugins/seasons.py
${PYSITELIB}/bravo/plugins/seasons.pyc
${PYSITELIB}/bravo/plugins/seasons.pyo
${PYSITELIB}/bravo/plugins/serializers.py
${PYSITELIB}/bravo/plugins/serializers.pyc
${PYSITELIB}/bravo/plugins/serializers.pyo
${PYSITELIB}/bravo/plugins/tracks.py
${PYSITELIB}/bravo/plugins/tracks.pyc
${PYSITELIB}/bravo/plugins/tracks.pyo
${PYSITELIB}/bravo/plugins/warp.py
${PYSITELIB}/bravo/plugins/warp.pyc
${PYSITELIB}/bravo/plugins/warp.pyo
${PYSITELIB}/bravo/plugins/web.py
${PYSITELIB}/bravo/plugins/web.pyc
${PYSITELIB}/bravo/plugins/web.pyo
${PYSITELIB}/bravo/policy/__init__.py
${PYSITELIB}/bravo/policy/__init__.pyc
${PYSITELIB}/bravo/policy/__init__.pyo
${PYSITELIB}/bravo/policy/dig.py
${PYSITELIB}/bravo/policy/dig.pyc
${PYSITELIB}/bravo/policy/dig.pyo
${PYSITELIB}/bravo/protocols/__init__.py
${PYSITELIB}/bravo/protocols/__init__.pyc
${PYSITELIB}/bravo/protocols/__init__.pyo
${PYSITELIB}/bravo/protocols/beta.py
${PYSITELIB}/bravo/protocols/beta.pyc
${PYSITELIB}/bravo/protocols/beta.pyo
${PYSITELIB}/bravo/protocols/infini.py
${PYSITELIB}/bravo/protocols/infini.pyc
${PYSITELIB}/bravo/protocols/infini.pyo
${PYSITELIB}/bravo/remote.py
${PYSITELIB}/bravo/remote.pyc
${PYSITELIB}/bravo/remote.pyo
${PYSITELIB}/bravo/service.py
${PYSITELIB}/bravo/service.pyc
${PYSITELIB}/bravo/service.pyo
${PYSITELIB}/bravo/simplex.py
${PYSITELIB}/bravo/simplex.pyc
${PYSITELIB}/bravo/simplex.pyo
${PYSITELIB}/bravo/stdio.py
${PYSITELIB}/bravo/stdio.pyc
${PYSITELIB}/bravo/stdio.pyo
${PYSITELIB}/bravo/terrain/__init__.py
${PYSITELIB}/bravo/terrain/__init__.pyc
${PYSITELIB}/bravo/terrain/__init__.pyo
${PYSITELIB}/bravo/terrain/trees.py
${PYSITELIB}/bravo/terrain/trees.pyc
${PYSITELIB}/bravo/terrain/trees.pyo
${PYSITELIB}/bravo/tests/__init__.py
${PYSITELIB}/bravo/tests/__init__.pyc
${PYSITELIB}/bravo/tests/__init__.pyo
${PYSITELIB}/bravo/tests/factories/__init__.py
${PYSITELIB}/bravo/tests/factories/__init__.pyc
${PYSITELIB}/bravo/tests/factories/__init__.pyo
${PYSITELIB}/bravo/tests/factories/test_beta.py
${PYSITELIB}/bravo/tests/factories/test_beta.pyc
${PYSITELIB}/bravo/tests/factories/test_beta.pyo
${PYSITELIB}/bravo/tests/plugins/__init__.py
${PYSITELIB}/bravo/tests/plugins/__init__.pyc
${PYSITELIB}/bravo/tests/plugins/__init__.pyo
${PYSITELIB}/bravo/tests/plugins/test_build_hooks.py
${PYSITELIB}/bravo/tests/plugins/test_build_hooks.pyc
${PYSITELIB}/bravo/tests/plugins/test_build_hooks.pyo
${PYSITELIB}/bravo/tests/plugins/test_commands.py
${PYSITELIB}/bravo/tests/plugins/test_commands.pyc
${PYSITELIB}/bravo/tests/plugins/test_commands.pyo
${PYSITELIB}/bravo/tests/plugins/test_generators.py
${PYSITELIB}/bravo/tests/plugins/test_generators.pyc
${PYSITELIB}/bravo/tests/plugins/test_generators.pyo
${PYSITELIB}/bravo/tests/plugins/test_physics.py
${PYSITELIB}/bravo/tests/plugins/test_physics.pyc
${PYSITELIB}/bravo/tests/plugins/test_physics.pyo
${PYSITELIB}/bravo/tests/plugins/test_recipes.py
${PYSITELIB}/bravo/tests/plugins/test_recipes.pyc
${PYSITELIB}/bravo/tests/plugins/test_recipes.pyo
${PYSITELIB}/bravo/tests/plugins/test_seasons.py
${PYSITELIB}/bravo/tests/plugins/test_seasons.pyc
${PYSITELIB}/bravo/tests/plugins/test_seasons.pyo
${PYSITELIB}/bravo/tests/plugins/test_serializers.py
${PYSITELIB}/bravo/tests/plugins/test_serializers.pyc
${PYSITELIB}/bravo/tests/plugins/test_serializers.pyo
${PYSITELIB}/bravo/tests/policy/__init__.py
${PYSITELIB}/bravo/tests/policy/__init__.pyc
${PYSITELIB}/bravo/tests/policy/__init__.pyo
${PYSITELIB}/bravo/tests/policy/test_dig.py
${PYSITELIB}/bravo/tests/policy/test_dig.pyc
${PYSITELIB}/bravo/tests/policy/test_dig.pyo
${PYSITELIB}/bravo/tests/protocols/__init__.py
${PYSITELIB}/bravo/tests/protocols/__init__.pyc
${PYSITELIB}/bravo/tests/protocols/__init__.pyo
${PYSITELIB}/bravo/tests/protocols/test_beta.py
${PYSITELIB}/bravo/tests/protocols/test_beta.pyc
${PYSITELIB}/bravo/tests/protocols/test_beta.pyo
${PYSITELIB}/bravo/tests/test_blocks.py
${PYSITELIB}/bravo/tests/test_blocks.pyc
${PYSITELIB}/bravo/tests/test_blocks.pyo
${PYSITELIB}/bravo/tests/test_chunk.py
${PYSITELIB}/bravo/tests/test_chunk.pyc
${PYSITELIB}/bravo/tests/test_chunk.pyo
${PYSITELIB}/bravo/tests/test_config.py
${PYSITELIB}/bravo/tests/test_config.pyc
${PYSITELIB}/bravo/tests/test_config.pyo
${PYSITELIB}/bravo/tests/test_entity.py
${PYSITELIB}/bravo/tests/test_entity.pyc
${PYSITELIB}/bravo/tests/test_entity.pyo
${PYSITELIB}/bravo/tests/test_inventory.py
${PYSITELIB}/bravo/tests/test_inventory.pyc
${PYSITELIB}/bravo/tests/test_inventory.pyo
${PYSITELIB}/bravo/tests/test_location.py
${PYSITELIB}/bravo/tests/test_location.pyc
${PYSITELIB}/bravo/tests/test_location.pyo
${PYSITELIB}/bravo/tests/test_nbt.py
${PYSITELIB}/bravo/tests/test_nbt.pyc
${PYSITELIB}/bravo/tests/test_nbt.pyo
${PYSITELIB}/bravo/tests/test_packets.py
${PYSITELIB}/bravo/tests/test_packets.pyc
${PYSITELIB}/bravo/tests/test_packets.pyo
${PYSITELIB}/bravo/tests/test_plugin.py
${PYSITELIB}/bravo/tests/test_plugin.pyc
${PYSITELIB}/bravo/tests/test_plugin.pyo
${PYSITELIB}/bravo/tests/test_simplex.py
${PYSITELIB}/bravo/tests/test_simplex.pyc
${PYSITELIB}/bravo/tests/test_simplex.pyo
${PYSITELIB}/bravo/tests/test_spatial.py
${PYSITELIB}/bravo/tests/test_spatial.pyc
${PYSITELIB}/bravo/tests/test_spatial.pyo
${PYSITELIB}/bravo/tests/test_utilities.py
${PYSITELIB}/bravo/tests/test_utilities.pyc
${PYSITELIB}/bravo/tests/test_utilities.pyo
${PYSITELIB}/bravo/tests/test_world.py
${PYSITELIB}/bravo/tests/test_world.pyc
${PYSITELIB}/bravo/tests/test_world.pyo
${PYSITELIB}/bravo/tests/utilities/__init__.py
${PYSITELIB}/bravo/tests/utilities/__init__.pyc
${PYSITELIB}/bravo/tests/utilities/__init__.pyo
${PYSITELIB}/bravo/tests/utilities/test_maths.py
${PYSITELIB}/bravo/tests/utilities/test_maths.pyc
${PYSITELIB}/bravo/tests/utilities/test_maths.pyo
${PYSITELIB}/bravo/utilities/__init__.py
${PYSITELIB}/bravo/utilities/__init__.pyc
${PYSITELIB}/bravo/utilities/__init__.pyo
${PYSITELIB}/bravo/utilities/bits.py
${PYSITELIB}/bravo/utilities/bits.pyc
${PYSITELIB}/bravo/utilities/bits.pyo
${PYSITELIB}/bravo/utilities/chat.py
${PYSITELIB}/bravo/utilities/chat.pyc
${PYSITELIB}/bravo/utilities/chat.pyo
${PYSITELIB}/bravo/utilities/coords.py
${PYSITELIB}/bravo/utilities/coords.pyc
${PYSITELIB}/bravo/utilities/coords.pyo
${PYSITELIB}/bravo/utilities/decos.py
${PYSITELIB}/bravo/utilities/decos.pyc
${PYSITELIB}/bravo/utilities/decos.pyo
${PYSITELIB}/bravo/utilities/maths.py
${PYSITELIB}/bravo/utilities/maths.pyc
${PYSITELIB}/bravo/utilities/maths.pyo
${PYSITELIB}/bravo/utilities/spatial.py
${PYSITELIB}/bravo/utilities/spatial.pyc
${PYSITELIB}/bravo/utilities/spatial.pyo
${PYSITELIB}/bravo/utilities/temporal.py
${PYSITELIB}/bravo/utilities/temporal.pyc
${PYSITELIB}/bravo/utilities/temporal.pyo
${PYSITELIB}/bravo/web.py
${PYSITELIB}/bravo/web.pyc
${PYSITELIB}/bravo/web.pyo
${PYSITELIB}/bravo/world.py
${PYSITELIB}/bravo/world.pyc
${PYSITELIB}/bravo/world.pyo
${PYSITELIB}/exocet/__init__.py
${PYSITELIB}/exocet/__init__.pyc
${PYSITELIB}/exocet/__init__.pyo
${PYSITELIB}/exocet/_components.py
${PYSITELIB}/exocet/_components.pyc
${PYSITELIB}/exocet/_components.pyo
${PYSITELIB}/exocet/_exocet.py
${PYSITELIB}/exocet/_exocet.pyc
${PYSITELIB}/exocet/_exocet.pyo
${PYSITELIB}/exocet/_filepath.py
${PYSITELIB}/exocet/_filepath.pyc
${PYSITELIB}/exocet/_filepath.pyo
${PYSITELIB}/exocet/_modules.py
${PYSITELIB}/exocet/_modules.pyc
${PYSITELIB}/exocet/_modules.pyo
${PYSITELIB}/exocet/_reflect.py
${PYSITELIB}/exocet/_reflect.pyc
${PYSITELIB}/exocet/_reflect.pyo
${PYSITELIB}/exocet/_runtime.py
${PYSITELIB}/exocet/_runtime.pyc
${PYSITELIB}/exocet/_runtime.pyo
${PYSITELIB}/exocet/_win32.py
${PYSITELIB}/exocet/_win32.pyc
${PYSITELIB}/exocet/_win32.pyo
${PYSITELIB}/exocet/_zippath.py
${PYSITELIB}/exocet/_zippath.pyc
${PYSITELIB}/exocet/_zippath.pyo
${PYSITELIB}/exocet/test/__init__.py
${PYSITELIB}/exocet/test/__init__.pyc
${PYSITELIB}/exocet/test/__init__.pyo
${PYSITELIB}/exocet/test/_ospathExample.py
${PYSITELIB}/exocet/test/_ospathExample.pyc
${PYSITELIB}/exocet/test/_ospathExample.pyo
${PYSITELIB}/exocet/test/test_exocet.py
${PYSITELIB}/exocet/test/test_exocet.pyc
${PYSITELIB}/exocet/test/test_exocet.pyo
${PYSITELIB}/exocet/test/test_modules.py
${PYSITELIB}/exocet/test/test_modules.pyc
${PYSITELIB}/exocet/test/test_modules.pyo
${PYSITELIB}/exocet/test/test_paths.py
${PYSITELIB}/exocet/test/test_paths.pyc
${PYSITELIB}/exocet/test/test_paths.pyo
${PYSITELIB}/exocet/test/testpackage/__init__.py
${PYSITELIB}/exocet/test/testpackage/__init__.pyc
${PYSITELIB}/exocet/test/testpackage/__init__.pyo
${PYSITELIB}/exocet/test/testpackage/baz.py
${PYSITELIB}/exocet/test/testpackage/baz.pyc
${PYSITELIB}/exocet/test/testpackage/baz.pyo
${PYSITELIB}/exocet/test/testpackage/foo.py
${PYSITELIB}/exocet/test/testpackage/foo.pyc
${PYSITELIB}/exocet/test/testpackage/foo.pyo
${PYSITELIB}/exocet/test/testpackage/topmodule.py
${PYSITELIB}/exocet/test/testpackage/topmodule.pyc
${PYSITELIB}/exocet/test/testpackage/topmodule.pyo
${PYSITELIB}/exocet/test/testpackage/util.py
${PYSITELIB}/exocet/test/testpackage/util.pyc
${PYSITELIB}/exocet/test/testpackage/util.pyo
${PYSITELIB}/exocet/test/testpackage_circular/__init__.py
${PYSITELIB}/exocet/test/testpackage_circular/__init__.pyc
${PYSITELIB}/exocet/test/testpackage_circular/__init__.pyo
${PYSITELIB}/exocet/test/testpackage_circular/baz.py
${PYSITELIB}/exocet/test/testpackage_circular/baz.pyc
${PYSITELIB}/exocet/test/testpackage_circular/baz.pyo
${PYSITELIB}/exocet/test/testpackage_circular/foo.py
${PYSITELIB}/exocet/test/testpackage_circular/foo.pyc
${PYSITELIB}/exocet/test/testpackage_circular/foo.pyo
${PYSITELIB}/exocet/test/testpackage_circular/topmodule.py
${PYSITELIB}/exocet/test/testpackage_circular/topmodule.pyc
${PYSITELIB}/exocet/test/testpackage_circular/topmodule.pyo
${PYSITELIB}/exocet/test/testpackage_circular/util.py
${PYSITELIB}/exocet/test/testpackage_circular/util.pyc
${PYSITELIB}/exocet/test/testpackage_circular/util.pyo
${PYSITELIB}/exocet/test/testpackage_localimports/__init__.py
${PYSITELIB}/exocet/test/testpackage_localimports/__init__.pyc
${PYSITELIB}/exocet/test/testpackage_localimports/__init__.pyo
${PYSITELIB}/exocet/test/testpackage_localimports/baz.py
${PYSITELIB}/exocet/test/testpackage_localimports/baz.pyc
${PYSITELIB}/exocet/test/testpackage_localimports/baz.pyo
${PYSITELIB}/exocet/test/testpackage_localimports/foo.py
${PYSITELIB}/exocet/test/testpackage_localimports/foo.pyc
${PYSITELIB}/exocet/test/testpackage_localimports/foo.pyo
${PYSITELIB}/exocet/test/testpackage_localimports/topmodule.py
${PYSITELIB}/exocet/test/testpackage_localimports/topmodule.pyc
${PYSITELIB}/exocet/test/testpackage_localimports/topmodule.pyo
${PYSITELIB}/exocet/test/testpackage_localimports/util.py
${PYSITELIB}/exocet/test/testpackage_localimports/util.pyc
${PYSITELIB}/exocet/test/testpackage_localimports/util.pyo
${PYSITELIB}/twisted/plugins/bravod.py
${PYSITELIB}/twisted/plugins/bravod.pyc
${PYSITELIB}/twisted/plugins/bravod.pyo
share/doc/bravo/administration.rst
share/doc/bravo/auxiliary.rst
share/doc/bravo/blocks.rst
share/doc/bravo/chunk.rst
share/doc/bravo/conf.py
share/doc/bravo/core.rst
share/doc/bravo/differences.rst
share/doc/bravo/entity.rst
share/doc/bravo/factory.rst
share/doc/bravo/features.rst
share/doc/bravo/globals.txt
share/doc/bravo/ibravo.rst
share/doc/bravo/index.rst
share/doc/bravo/location.rst
share/doc/bravo/noise.rst
share/doc/bravo/packets.rst
share/doc/bravo/philosophy.rst
share/doc/bravo/plugin.rst
share/doc/bravo/plugins.rst
share/doc/bravo/protocol.rst
share/doc/bravo/simplex.rst
share/doc/bravo/stdio.rst
share/doc/bravo/tools.rst
share/doc/bravo/utilities.rst
share/doc/bravo/web.rst
share/doc/bravo/world.rst
share/examples/bravo/bravo.ini
share/examples/rc.d/bravo
@pkgdir etc/bravo
