# uptimed/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2007/02/27 22:54:17 dzoe Exp $
bin/uprecords
lib/libuptimed.la
man/man1/uprecords.1
man/man8/uptimed.8
sbin/uptimed
share/examples/rc.d/uptimed
