# xtreemfs-server/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2011/11/08 01:42:14 tnn2 Exp $
bin/xtfs_chstatus
bin/xtfs_cleanup
bin/xtfs_mrcdbtool
bin/xtfs_remove_osd
bin/xtfs_scrub
bin/xtfs_snap
man/man1/xtfs_chstatus.1
man/man1/xtfs_cleanup.1
man/man1/xtfs_mrcdbtool.1
man/man1/xtfs_remove_osd.1
man/man1/xtfs_scrub.1
man/man1/xtfs_snap.1
share/doc/xtreemfs-server/LICENSE
share/doc/xtreemfs-tools/LICENSE
share/examples/rc.d/xtreemfs_dir
share/examples/rc.d/xtreemfs_mrc
share/examples/rc.d/xtreemfs_osd
share/examples/xtreemfs/dirconfig.properties
share/examples/xtreemfs/mrcconfig.properties
share/examples/xtreemfs/osdconfig.properties
share/examples/xtreemfs/server-repl-plugin/dir.properties
share/examples/xtreemfs/server-repl-plugin/mrc.properties
share/xtreemfs/java/lib/BabuDB.jar
share/xtreemfs/java/lib/Flease.jar
share/xtreemfs/java/lib/Foundation.jar
share/xtreemfs/java/lib/XtreemFS.jar
share/xtreemfs/java/lib/jdmkrt.jar
share/xtreemfs/java/lib/jdmktk.jar
share/xtreemfs/java/lib/protobuf-java-2.3.0.jar
share/xtreemfs/server-repl-plugin/lib/Flease.jar
share/xtreemfs/server-repl-plugin/lib/PBRPC.jar
share/xtreemfs/server-repl-plugin/lib/protobuf-java-2.3.0.jar
share/xtreemfs/server-repl-plugin/replication.jar
