# ganglia-monitor-core319/PLIST.Linux
share/examples/rc.d/gmetad.init
share/examples/rc.d/gmond.init
