# drweb/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2003/04/10 18:51:52 jmmv Exp $
bin/drweb
bin/drweb-fbsd
bin/drweb-update
etc/rc.d/drwebd
sbin/drwebd
sbin/drwebd-fbsd
${_DOCDIR}/ChangeLog
${_DOCDIR}/README
${_DOCDIR}/drwebd-api
${_DOCDIR}/readme.daemon
${_DOCDIR}/readme.eicar
${_DOCDIR}/readme.scanner
${_DOCDIR}/readme.update
${_EGDIR}/drweb.key
${_EGDIR}/drweb32.dll
${_EGDIR}/drweb32.ini
${_EGDIR}/drwebase.vdb
${_EGDIR}/drwebd.key
@dirrm ${_DOCDIR}
@dirrm ${_EGDIR}
