# drweb/PLIST
@comment $NetBSD: PLIST,v 1.4 2003/09/16 16:57:59 bubuchka Exp $
bin/drweb
bin/drweb-bin
etc/rc.d/drwebd
sbin/drweb-update
sbin/drwebd
sbin/drwebd-bin
${DRW_DOCDIR}/ChangeLog
${DRW_DOCDIR}/README
${DRW_DOCDIR}/drwebd-api
${DRW_DOCDIR}/readme.daemon
${DRW_DOCDIR}/readme.eicar
${DRW_DOCDIR}/readme.scanner
${DRW_DOCDIR}/readme.update
${DRW_EGDIR}/drweb.key
${DRW_EGDIR}/drweb32.dll
${DRW_EGDIR}/drweb32.ini
${DRW_EGDIR}/drwebase.vdb
${DRW_EGDIR}/drwebd.key
@dirrm ${DRW_DOCDIR}
@dirrm ${DRW_EGDIR}
