# dbmail/PLIST
@comment $NetBSD: PLIST,v 1.5 2013/09/06 15:11:01 tnn2 Exp $
lib/dbmail/libauth_ldap.la
lib/dbmail/libauth_sql.la
lib/dbmail/libdbmail.la
${PLIST.sieve}lib/dbmail/libsort_sieve.la
man/man1/dbmail-deliver.1
man/man1/dbmail.1
man/man5/dbmail.conf.5
man/man8/dbmail-export.8
man/man8/dbmail-httpd.8
man/man8/dbmail-imapd.8
man/man8/dbmail-lmtpd.8
man/man8/dbmail-pop3d.8
man/man8/dbmail-sievecmd.8
man/man8/dbmail-timsieved.8
man/man8/dbmail-users.8
man/man8/dbmail-util.8
sbin/dbmail-deliver
sbin/dbmail-export
sbin/dbmail-httpd
sbin/dbmail-imapd
sbin/dbmail-lmtpd
sbin/dbmail-pop3d
${PLIST.sieve}sbin/dbmail-sievecmd
${PLIST.sieve}sbin/dbmail-timsieved
sbin/dbmail-users
sbin/dbmail-util
share/dbmail/sql/mysql/2_1_7-2_2_0.mysql
share/dbmail/sql/mysql/2_2-3_0.mysql
share/dbmail/sql/mysql/2_2_5-2_3_0.mysql
share/dbmail/sql/mysql/2_3_0-2_3_2.mysql
share/dbmail/sql/mysql/2_3_4-2_3_5.mysql
share/dbmail/sql/mysql/2_3_5-2_3_6.mysql
share/dbmail/sql/mysql/2_3_6-3_0_0.mysql
share/dbmail/sql/mysql/3_0_2-3_1_0.mysql
share/dbmail/sql/mysql/create_tables.mysql
share/dbmail/sql/mysql/fix_foreign_keys.mysql
share/dbmail/sql/mysql/migrate_from_1.x_to_2.0_innodb.mysql
share/dbmail/sql/mysql/migrate_from_2.0_to_2.2.mysql
share/dbmail/sql/mysql/update_headervalue_01.mysql
share/dbmail/sql/postgresql/2_1_7-2_2_0.pgsql
share/dbmail/sql/postgresql/2_2-3_0.pgsql
share/dbmail/sql/postgresql/2_2_5-2_3_0.pgsql
share/dbmail/sql/postgresql/2_3_4-2_3_5.pgsql
share/dbmail/sql/postgresql/2_3_5-2_3_6.pgsql
share/dbmail/sql/postgresql/2_3_6-3_0_0.pgsql
share/dbmail/sql/postgresql/3_0_2-3_1_0.psql
share/dbmail/sql/postgresql/create_tables.pgsql
share/dbmail/sql/postgresql/migrate_from_1.x_to_2.0.pgsql
share/dbmail/sql/postgresql/migrate_from_2.0_to_2.2.pgsql
share/dbmail/sql/sqlite/2_1_7-2_2_0.sqlite
share/dbmail/sql/sqlite/2_2-3_0.sqlite
share/dbmail/sql/sqlite/2_2_5-2_3_0.sqlite
share/dbmail/sql/sqlite/2_3_4-2_3_5.sqlite
share/dbmail/sql/sqlite/2_3_5-2_3_6.sqlite
share/dbmail/sql/sqlite/2_3_6-3_0_0.sqlite
share/dbmail/sql/sqlite/3_0_2-3_1_0.sqlite
share/dbmail/sql/sqlite/create_tables.sqlite
share/dbmail/sql/sqlite/trigger.tmpl.sql
share/doc/dbmail/AUTHORS
share/doc/dbmail/BUGS
share/doc/dbmail/COPYING
share/doc/dbmail/INSTALL
share/doc/dbmail/NEWS
share/doc/dbmail/README
share/doc/dbmail/README.aliases
share/doc/dbmail/README.exim
share/doc/dbmail/README.filters
share/doc/dbmail/README.http
share/doc/dbmail/README.ldap
share/doc/dbmail/README.postfix
share/doc/dbmail/README.qmail
share/doc/dbmail/README.sieve
share/doc/dbmail/README.smtp
share/doc/dbmail/README.solaris
share/doc/dbmail/README.usermap
share/doc/dbmail/THANKS
share/doc/dbmail/UPGRADING
share/examples/dbmail/dbmail.conf
share/examples/rc.d/dbmail
share/examples/rc.d/dbmailimapd
share/examples/rc.d/dbmaillmtpd
share/examples/rc.d/dbmailpop3d
${PLIST.sieve}share/examples/rc.d/dbmailtimsieved
