# nginx/PLIST
@comment $NetBSD: PLIST,v 1.6 2008/05/06 11:32:22 rootcho Exp $
sbin/nginx
share/examples/nginx/conf/fastcgi_params
share/examples/nginx/conf/koi-utf
share/examples/nginx/conf/koi-win
share/examples/nginx/conf/mime.types
share/examples/nginx/conf/nginx.conf
share/examples/nginx/conf/win-utf
share/examples/nginx/html/50x.html
share/examples/nginx/html/index.html
share/examples/rc.d/nginx
@dirrm share/examples/nginx/html
@dirrm share/examples/nginx/conf
@dirrm share/examples/nginx
