# linux-nfs-utils/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:00 jeremy-c-reed Exp $
man/man5/exports.5
man/man7/nfsd.7
man/man8/exportfs.8
man/man8/lockd.8
man/man8/mountd.8
man/man8/nfsd.8
man/man8/nfsstat.8
man/man8/nhfsgraph.8
man/man8/nhfsnums.8
man/man8/nhfsrun.8
man/man8/nhfsstone.8
man/man8/rpc.lockd.8
man/man8/rpc.mountd.8
man/man8/rpc.nfsd.8
man/man8/rpc.rquotad.8
man/man8/rpc.statd.8
man/man8/rquotad.8
man/man8/showmount.8
man/man8/statd.8
sbin/exportfs
sbin/nfsstat
sbin/nhfsgraph
sbin/nhfsnums
sbin/nhfsrun
sbin/nhfsstone
sbin/rpc.lockd
sbin/rpc.mountd
sbin/rpc.nfsd
sbin/rpc.rquotad
sbin/rpc.statd
sbin/showmount
share/examples/rc.d/mountd
share/examples/rc.d/nfsd
share/examples/rc.d/nfslocking
