# py-jabber-proxy65/PLIST
@comment $NetBSD: PLIST,v 1.3 2008/08/13 04:56:31 schnoebe Exp $
${PYSITELIB}/proxy65/__init__.py
${PYSITELIB}/proxy65/__init__.pyc
${PYSITELIB}/proxy65/__init__.pyo
${PYSITELIB}/proxy65/proxy65.py
${PYSITELIB}/proxy65/proxy65.pyc
${PYSITELIB}/proxy65/proxy65.pyo
${PYSITELIB}/proxy65/socks5.py
${PYSITELIB}/proxy65/socks5.pyc
${PYSITELIB}/proxy65/socks5.pyo
${PYSITELIB}/proxy65/tap.py
${PYSITELIB}/proxy65/tap.pyc
${PYSITELIB}/proxy65/tap.pyo
${PYSITELIB}/twisted/plugins/proxy65.py
${PYSITELIB}/twisted/plugins/proxy65.pyc
${PYSITELIB}/twisted/plugins/proxy65.pyo
share/doc/proxy65/README
share/examples/rc.d/proxy65
@dirrm share/doc/proxy65
@dirrm ${PYSITELIB}/proxy65
