# py-jabber-idavoll/PLIST
@comment $NetBSD: PLIST,v 1.1 2013/12/01 21:09:50 schnoebe Exp $
${PYSITELIB}/${EGG_INFODIR}/PKG-INFO
${PYSITELIB}/${EGG_INFODIR}/SOURCES.txt
${PYSITELIB}/${EGG_INFODIR}/dependency_links.txt
${PYSITELIB}/${EGG_INFODIR}/not-zip-safe
${PYSITELIB}/${EGG_INFODIR}/requires.txt
${PYSITELIB}/${EGG_INFODIR}/top_level.txt
${PYSITELIB}/idavoll/__init__.py
${PYSITELIB}/idavoll/__init__.pyc
${PYSITELIB}/idavoll/__init__.pyo
${PYSITELIB}/idavoll/backend.py
${PYSITELIB}/idavoll/backend.pyc
${PYSITELIB}/idavoll/backend.pyo
${PYSITELIB}/idavoll/error.py
${PYSITELIB}/idavoll/error.pyc
${PYSITELIB}/idavoll/error.pyo
${PYSITELIB}/idavoll/gateway.py
${PYSITELIB}/idavoll/gateway.pyc
${PYSITELIB}/idavoll/gateway.pyo
${PYSITELIB}/idavoll/iidavoll.py
${PYSITELIB}/idavoll/iidavoll.pyc
${PYSITELIB}/idavoll/iidavoll.pyo
${PYSITELIB}/idavoll/memory_storage.py
${PYSITELIB}/idavoll/memory_storage.pyc
${PYSITELIB}/idavoll/memory_storage.pyo
${PYSITELIB}/idavoll/pgsql_storage.py
${PYSITELIB}/idavoll/pgsql_storage.pyc
${PYSITELIB}/idavoll/pgsql_storage.pyo
${PYSITELIB}/idavoll/tap.py
${PYSITELIB}/idavoll/tap.pyc
${PYSITELIB}/idavoll/tap.pyo
${PYSITELIB}/idavoll/tap_http.py
${PYSITELIB}/idavoll/tap_http.pyc
${PYSITELIB}/idavoll/tap_http.pyo
${PYSITELIB}/idavoll/test/__init__.py
${PYSITELIB}/idavoll/test/__init__.pyc
${PYSITELIB}/idavoll/test/__init__.pyo
${PYSITELIB}/idavoll/test/test_backend.py
${PYSITELIB}/idavoll/test/test_backend.pyc
${PYSITELIB}/idavoll/test/test_backend.pyo
${PYSITELIB}/idavoll/test/test_gateway.py
${PYSITELIB}/idavoll/test/test_gateway.pyc
${PYSITELIB}/idavoll/test/test_gateway.pyo
${PYSITELIB}/idavoll/test/test_storage.py
${PYSITELIB}/idavoll/test/test_storage.pyc
${PYSITELIB}/idavoll/test/test_storage.pyo
${PYSITELIB}/twisted/plugins/idavoll.py
${PYSITELIB}/twisted/plugins/idavoll.pyc
${PYSITELIB}/twisted/plugins/idavoll.pyo
${PYSITELIB}/twisted/plugins/idavoll_http.py
${PYSITELIB}/twisted/plugins/idavoll_http.pyc
${PYSITELIB}/twisted/plugins/idavoll_http.pyo
share/doc/idavoll/INSTALL
share/doc/idavoll/LICENSE
share/doc/idavoll/NEWS
share/doc/idavoll/README
share/examples/rc.d/idavoll
share/idavoll/gateway.sql
share/idavoll/pubsub.sql
share/idavoll/to_idavoll_0.8.sql
@pkgdir share/examples/jabber
