# noip/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/10/11 10:45:00 thomasklausner Exp $
bin/noip2
share/doc/noip/COPYING
share/doc/noip/README.FIRST
share/examples/rc.d/noip
