# slim/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/10/11 10:45:18 thomasklausner Exp $
sbin/slim
man/man1/slim.1
share/examples/rc.d/slim
share/examples/slim.conf
share/slim/themes/default/background.jpg
share/slim/themes/default/panel.png
share/slim/themes/default/slim.theme
