# libreswan-5/PLIST
@comment $NetBSD$
libexec/ipsec/_asn1check
libexec/ipsec/_dncheck
libexec/ipsec/_enumcheck
libexec/ipsec/_hunkcheck
libexec/ipsec/_import_crl
libexec/ipsec/_ipcheck
libexec/ipsec/_jambufcheck
libexec/ipsec/_keyidcheck
libexec/ipsec/_plutorun
libexec/ipsec/_timecheck
libexec/ipsec/_ttodatacheck
libexec/ipsec/_unbound-hook
libexec/ipsec/_updown
libexec/ipsec/_updown.bsd
libexec/ipsec/_vendoridcheck
libexec/ipsec/addconn
libexec/ipsec/algparse
libexec/ipsec/cavp
libexec/ipsec/ecdsasigkey
libexec/ipsec/letsencrypt
libexec/ipsec/newhostkey
libexec/ipsec/pluto
libexec/ipsec/readwriteconf
libexec/ipsec/rsasigkey
libexec/ipsec/setup
libexec/ipsec/showhostkey
libexec/ipsec/showroute
libexec/ipsec/whack
man/man5/ipsec.conf.5
man/man5/ipsec.secrets.5
man/man7/libreswan.7
man/man8/ipsec-_plutorun.8
man/man8/ipsec-_updown.8
man/man8/ipsec-_updown.bsd.8
man/man8/ipsec-add.8
man/man8/ipsec-addconn.8
man/man8/ipsec-algparse.8
man/man8/ipsec-briefconnectionstatus.8
man/man8/ipsec-briefstatus.8
man/man8/ipsec-certutil.8
man/man8/ipsec-checkconfig.8
man/man8/ipsec-checknflog.8
man/man8/ipsec-checknss.8
man/man8/ipsec-connectionstatus.8
man/man8/ipsec-crlutil.8
man/man8/ipsec-delete.8
man/man8/ipsec-down.8
man/man8/ipsec-ecdsasigkey.8
man/man8/ipsec-fetchcrls.8
man/man8/ipsec-fipsstatus.8
man/man8/ipsec-globalstatus.8
man/man8/ipsec-import.8
man/man8/ipsec-initnss.8
man/man8/ipsec-letsencrypt.8
man/man8/ipsec-listall.8
man/man8/ipsec-listcacerts.8
man/man8/ipsec-listcerts.8
man/man8/ipsec-listcrls.8
man/man8/ipsec-listen.8
man/man8/ipsec-listpubkeys.8
man/man8/ipsec-modutil.8
man/man8/ipsec-newhostkey.8
man/man8/ipsec-ondemand.8
man/man8/ipsec-pk12util.8
man/man8/ipsec-pluto.8
man/man8/ipsec-purgeocsp.8
man/man8/ipsec-readwriteconf.8
man/man8/ipsec-redirect.8
man/man8/ipsec-replace.8
man/man8/ipsec-rereadall.8
man/man8/ipsec-rereadcerts.8
man/man8/ipsec-rereadsecrets.8
man/man8/ipsec-restart.8
man/man8/ipsec-route.8
man/man8/ipsec-rsasigkey.8
man/man8/ipsec-setup.8
man/man8/ipsec-showhostkey.8
man/man8/ipsec-showroute.8
man/man8/ipsec-showstates.8
man/man8/ipsec-shuntstatus.8
man/man8/ipsec-start.8
man/man8/ipsec-status.8
man/man8/ipsec-stop.8
man/man8/ipsec-trafficstatus.8
man/man8/ipsec-unroute.8
man/man8/ipsec-up.8
man/man8/ipsec-vfychain.8
man/man8/ipsec-whack.8
man/man8/ipsec.8
man/man8/pluto.8
sbin/ipsec
share/examples/libreswan/ipsec.conf-sample
share/examples/libreswan/ipsec.d/policies/block
share/examples/libreswan/ipsec.d/policies/clear
share/examples/libreswan/ipsec.d/policies/clear-or-private
share/examples/libreswan/ipsec.d/policies/portexcludes.conf
share/examples/libreswan/ipsec.d/policies/private
share/examples/libreswan/ipsec.d/policies/private-or-clear
share/examples/libreswan/ipsec.secrets-sample
share/examples/libreswan/logrotate.d/libreswan
share/examples/libreswan/pam.d/pluto
share/examples/libreswan/rc.d/pluto.sh
