# alphamail/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2008/02/12 11:36:29 shattered Exp $
lib/alphamail/AlphaMail.pm
lib/alphamail/AlphaMail/AddressBook.pm
lib/alphamail/AlphaMail/AdminAccessRule.pm
lib/alphamail/AlphaMail/AdminUsers.pm
lib/alphamail/AlphaMail/Config.pm
lib/alphamail/AlphaMail/Filter.pm
lib/alphamail/AlphaMail/HTMLHelper.pm
lib/alphamail/AlphaMail/MailTool.pm
lib/alphamail/AlphaMail/Message.pm
lib/alphamail/AlphaMail/Middleware.pm
lib/alphamail/AlphaMail/Preferences.pm
lib/alphamail/AlphaMail/Quota.pm
lib/alphamail/AlphaMail/RequestHandler.pm
lib/alphamail/AlphaMail/SAPMConfig.pm
lib/alphamail/AlphaMail/SessionUtil.pm
lib/alphamail/AlphaMail/Signature.pm
lib/alphamail/AlphaMail/SpamConfig.pm
lib/alphamail/Apache2/SiteControl.pm
lib/alphamail/Apache2/SiteControl/GrantAllRule.pm
lib/alphamail/Apache2/SiteControl/ManagerFactory.pm
lib/alphamail/Apache2/SiteControl/PermissionManager.pm
lib/alphamail/Apache2/SiteControl/Rule.pm
lib/alphamail/Apache2/SiteControl/User.pm
lib/alphamail/Apache2/SiteControl/UserFactory.pm
lib/alphamail/Encode/IMAPUTF7.pm
lib/alphamail/IMAPAuth.pm
lib/alphamail/ManagerFactory.pm
libexec/sandbox
sbin/alphamail_genconfig
sbin/alphamail_stats
sbin/garbage_sweeper
sbin/imap_webcache
share/alphamail/conf/README
share/alphamail/conf/alphamail.conf
share/alphamail/conf/init.pl
share/alphamail/conf/log4perl.conf
share/alphamail/docs/About
share/alphamail/docs/AlphamailAdmin.pod
share/alphamail/docs/Dependencies
share/alphamail/docs/build.sh
share/alphamail/docs/caching_algorithm
share/alphamail/docs/certificate_authorities.pem
share/alphamail/docs/middleware_protocol.txt
share/alphamail/docs/notes/OSX
share/alphamail/docs/notes/SunStudio
share/alphamail/docs/notes/debugging
share/alphamail/htdocs/download.html
share/alphamail/htdocs/error.html
share/alphamail/htdocs/filevault.html
share/alphamail/htdocs/index.html
share/alphamail/htdocs/mail/address_mail.html
share/alphamail/htdocs/mail/addressbook/edit.mhtml
share/alphamail/htdocs/mail/addressbook/export_csv.html
share/alphamail/htdocs/mail/addressbook/import.mhtml
share/alphamail/htdocs/mail/addressbook/import_complete.html
share/alphamail/htdocs/mail/addressbook/import_csv.html
share/alphamail/htdocs/mail/addressbook/import_csv_save.html
share/alphamail/htdocs/mail/addressbook/import_fname.html
share/alphamail/htdocs/mail/addressbook/import_imho.html
share/alphamail/htdocs/mail/addressbook/import_lname.html
share/alphamail/htdocs/mail/addressbook/index.html
share/alphamail/htdocs/mail/addressbook/lists.mhtml
share/alphamail/htdocs/mail/addressbook/process_edit.html
share/alphamail/htdocs/mail/addressbook/process_lists.html
share/alphamail/htdocs/mail/addressbook/process_take.html
share/alphamail/htdocs/mail/addressbook/take.html
share/alphamail/htdocs/mail/addresslist.html
share/alphamail/htdocs/mail/admin/autohandler
share/alphamail/htdocs/mail/admin/guide.html
share/alphamail/htdocs/mail/admin/index.html
share/alphamail/htdocs/mail/admin/quotas.jpg
share/alphamail/htdocs/mail/admin/show_stats.mhtml
share/alphamail/htdocs/mail/admin/test.mhtml
share/alphamail/htdocs/mail/autohandler
share/alphamail/htdocs/mail/check_spelling.html
share/alphamail/htdocs/mail/compose.html
share/alphamail/htdocs/mail/docs/faq/attach1.gif
share/alphamail/htdocs/mail/docs/faq/attach2.gif
share/alphamail/htdocs/mail/docs/faq/attach3.gif
share/alphamail/htdocs/mail/docs/faq/colors.gif
share/alphamail/htdocs/mail/docs/faq/dates1.gif
share/alphamail/htdocs/mail/docs/faq/drafts.gif
share/alphamail/htdocs/mail/docs/faq/filter.gif
share/alphamail/htdocs/mail/docs/faq/index.html
share/alphamail/htdocs/mail/docs/faq/notrash1.gif
share/alphamail/htdocs/mail/docs/faq/notrash2.gif
share/alphamail/htdocs/mail/docs/faq/notrash3.gif
share/alphamail/htdocs/mail/docs/faq/other.html
share/alphamail/htdocs/mail/docs/faq/sort.gif
share/alphamail/htdocs/mail/docs/faq/timeout.gif
share/alphamail/htdocs/mail/docs/index.html
share/alphamail/htdocs/mail/docs/new_features.html
share/alphamail/htdocs/mail/feedback.html
share/alphamail/htdocs/mail/first_login.html
share/alphamail/htdocs/mail/fix_spelling.html
share/alphamail/htdocs/mail/folderlist.mhtml
share/alphamail/htdocs/mail/footer.mhtml
share/alphamail/htdocs/mail/get_attachment.html
share/alphamail/htdocs/mail/header.mhtml
share/alphamail/htdocs/mail/help.html
share/alphamail/htdocs/mail/index.html
share/alphamail/htdocs/mail/logout.html
share/alphamail/htdocs/mail/menu.mhtml
share/alphamail/htdocs/mail/other_folders.html
share/alphamail/htdocs/mail/process_compose.html
share/alphamail/htdocs/mail/process_first_login.html
share/alphamail/htdocs/mail/process_messages.html
share/alphamail/htdocs/mail/process_read_message.html
share/alphamail/htdocs/mail/process_search.html
share/alphamail/htdocs/mail/quota.html
share/alphamail/htdocs/mail/quota_graph.html
share/alphamail/htdocs/mail/raw.html
share/alphamail/htdocs/mail/read.html
share/alphamail/htdocs/mail/remove_attachment.html
share/alphamail/htdocs/mail/renew_session.html
share/alphamail/htdocs/mail/resume_compose.html
share/alphamail/htdocs/mail/sanitized.html
share/alphamail/htdocs/mail/search.html
share/alphamail/htdocs/mail/search_results.html
share/alphamail/htdocs/mail/settings/add_filter.html
share/alphamail/htdocs/mail/settings/altfolder_update.html
share/alphamail/htdocs/mail/settings/altprefix_update.html
share/alphamail/htdocs/mail/settings/edit_folders.html
share/alphamail/htdocs/mail/settings/export.mhtml
share/alphamail/htdocs/mail/settings/export_mail.html
share/alphamail/htdocs/mail/settings/export_verify.html
share/alphamail/htdocs/mail/settings/filters.mhtml
share/alphamail/htdocs/mail/settings/folders.mhtml
share/alphamail/htdocs/mail/settings/general.mhtml
share/alphamail/htdocs/mail/settings/index.html
share/alphamail/htdocs/mail/settings/share.mhtml
share/alphamail/htdocs/mail/settings/sharing_agree.html
share/alphamail/htdocs/mail/settings/signatures.mhtml
share/alphamail/htdocs/mail/settings/spam.mhtml
share/alphamail/htdocs/mail/settings/update_filters.html
share/alphamail/htdocs/mail/settings/update_folderlist.html
share/alphamail/htdocs/mail/settings/update_general.html
share/alphamail/htdocs/mail/settings/update_signatures.html
share/alphamail/htdocs/mail/settings/update_spam.html
share/alphamail/htdocs/mail/share_options.html
share/alphamail/htdocs/mail/share_upload.html
share/alphamail/htdocs/mail/verify_password.html
share/alphamail/htdocs/mail/view_attachment.html
share/alphamail/htdocs/mail/viewers/excel.html
share/alphamail/htdocs/mail/viewers/html.html
share/alphamail/htdocs/mail/viewers/targz.html
share/alphamail/htdocs/mail/viewers/text.html
share/alphamail/htdocs/mail/viewers/word.html
share/alphamail/htdocs/mail/viewers/zip.html
share/alphamail/htdocs/static/amail_icon.jpg
share/alphamail/htdocs/static/amail_large.jpg
share/alphamail/htdocs/static/bell.wav
share/alphamail/htdocs/static/default.css
share/alphamail/htdocs/static/paperclip.png
share/alphamail/htdocs/static/rte/blank.htm
share/alphamail/htdocs/static/rte/html2xhtml.js
share/alphamail/htdocs/static/rte/images/bgcolor.gif
share/alphamail/htdocs/static/rte/images/blackdot.gif
share/alphamail/htdocs/static/rte/images/bold.gif
share/alphamail/htdocs/static/rte/images/centre.gif
share/alphamail/htdocs/static/rte/images/copy.gif
share/alphamail/htdocs/static/rte/images/cut.gif
share/alphamail/htdocs/static/rte/images/email.gif
share/alphamail/htdocs/static/rte/images/hr.gif
share/alphamail/htdocs/static/rte/images/hyperlink.gif
share/alphamail/htdocs/static/rte/images/image.gif
share/alphamail/htdocs/static/rte/images/indent.gif
share/alphamail/htdocs/static/rte/images/insert_table.gif
share/alphamail/htdocs/static/rte/images/italic.gif
share/alphamail/htdocs/static/rte/images/justifyfull.gif
share/alphamail/htdocs/static/rte/images/left_just.gif
share/alphamail/htdocs/static/rte/images/list.gif
share/alphamail/htdocs/static/rte/images/numbered_list.gif
share/alphamail/htdocs/static/rte/images/outdent.gif
share/alphamail/htdocs/static/rte/images/paste.gif
share/alphamail/htdocs/static/rte/images/redo.gif
share/alphamail/htdocs/static/rte/images/right_just.gif
share/alphamail/htdocs/static/rte/images/smiley.gif
share/alphamail/htdocs/static/rte/images/spellcheck.gif
share/alphamail/htdocs/static/rte/images/textcolor.gif
share/alphamail/htdocs/static/rte/images/underline.gif
share/alphamail/htdocs/static/rte/images/undo.gif
share/alphamail/htdocs/static/rte/insert_link.htm
share/alphamail/htdocs/static/rte/insert_table.htm
share/alphamail/htdocs/static/rte/multi.htm
share/alphamail/htdocs/static/rte/palette.htm
share/alphamail/htdocs/static/rte/richtext.js
share/alphamail/htdocs/static/rte/richtext_compressed.js
share/alphamail/htdocs/static/rte/rte.css
share/alphamail/htdocs/static/wick.css
share/alphamail/htdocs/static/wick.js
share/alphamail/htdocs/static/wz_tooltip.js
share/alphamail/tests/Makefile
share/alphamail/tests/README
share/alphamail/tests/cache/append
share/alphamail/tests/cache/clear_inbox
share/alphamail/tests/cache/filters
share/alphamail/tests/cache/folders
share/alphamail/tests/cache/messagelists
share/alphamail/tests/cache/messagemoves
share/alphamail/tests/cache/search
share/alphamail/tests/cache/setup
share/alphamail/tests/cache/test
share/alphamail/tests/interface/addlist
share/alphamail/tests/interface/attachments
share/alphamail/tests/interface/autobcc
share/alphamail/tests/interface/compose
share/alphamail/tests/interface/contacts
share/alphamail/tests/interface/drafts
share/alphamail/tests/interface/faq.txt
share/alphamail/tests/interface/first_login
share/alphamail/tests/interface/flagstest
share/alphamail/tests/interface/folders
share/alphamail/tests/interface/forwardtest
share/alphamail/tests/interface/helplink
share/alphamail/tests/interface/index
share/alphamail/tests/interface/init
share/alphamail/tests/interface/largefile
share/alphamail/tests/interface/leftmenu
share/alphamail/tests/interface/linkify
share/alphamail/tests/interface/login
share/alphamail/tests/interface/logout
share/alphamail/tests/interface/multiplesessions
share/alphamail/tests/interface/njs_login
share/alphamail/tests/interface/readoptions
share/alphamail/tests/interface/sendtest
share/alphamail/tests/interface/sorting
share/alphamail/tests/interface/soundnotification
share/alphamail/tests/interface/spellcheck
share/alphamail/tests/interface/test
share/alphamail/tests/interface/util.pm
share/alphamail/tests/lib/TestUser.pm
share/alphamail/tests/lib/address_book
share/alphamail/tests/lib/mail_tools
share/alphamail/tests/lib/message/base
share/alphamail/tests/lib/message/resources/html
share/alphamail/tests/lib/message/resources/iso2022-jp
share/alphamail/tests/lib/message/resources/tarfile
share/alphamail/tests/lib/message/resources/utf8-body
share/alphamail/tests/lib/message/resources/utf8-jp
share/alphamail/tests/lib/message/resources/utf8-jp-img
share/alphamail/tests/lib/message/runtests
share/alphamail/tests/lib/message/send_fails
share/alphamail/tests/lib/signatures
share/alphamail/tests/lib/spam_config
share/alphamail/tests/lib/test
share/alphamail/tests/lib/urlbuilder
share/alphamail/tests/resources/hurricane.gif
share/alphamail/tests/resources/hurricane.jpg
share/alphamail/tests/resources/lightning.gif
share/alphamail/tests/resources/lightning.jpg
share/alphamail/tests/resources/lightning.png
share/alphamail/tests/resources/lightning.tif
share/alphamail/tests/resources/splash.gif
share/alphamail/tests/resources/splash.jpg
share/alphamail/tests/test
share/alphamail/util/hang_detector
share/alphamail/util/init.d/alphamail_tmpfs
share/alphamail/util/init.d/imap_webcache
share/examples/rc.d/imap_webcache
@dirrm share/alphamail/util/init.d
@dirrm share/alphamail/util
@dirrm share/alphamail/tests/resources
@dirrm share/alphamail/tests/lib/message/resources
@dirrm share/alphamail/tests/lib/message
@dirrm share/alphamail/tests/lib
@dirrm share/alphamail/tests/interface
@dirrm share/alphamail/tests/cache
@dirrm share/alphamail/tests
@dirrm share/alphamail/htdocs/static/rte/images
@dirrm share/alphamail/htdocs/static/rte
@dirrm share/alphamail/htdocs/static
@dirrm share/alphamail/htdocs/mail/viewers
@dirrm share/alphamail/htdocs/mail/settings
@dirrm share/alphamail/htdocs/mail/docs/faq
@dirrm share/alphamail/htdocs/mail/docs
@dirrm share/alphamail/htdocs/mail/admin
@dirrm share/alphamail/htdocs/mail/addressbook
@dirrm share/alphamail/htdocs/mail
@dirrm share/alphamail/htdocs
@dirrm share/alphamail/docs/notes
@dirrm share/alphamail/docs
@dirrm share/alphamail/conf
@dirrm lib/alphamail/Encode
@dirrm lib/alphamail/Apache2/SiteControl
@dirrm lib/alphamail/Apache2
@dirrm lib/alphamail/AlphaMail
@dirrm lib/alphamail
