# ypserv/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:08 jeremy-c-reed Exp $
include/rpcsvc/ypxfrd.x
sbin/create_printcap
sbin/makedbm
sbin/match_printcap
sbin/mknetid
sbin/pwupdate
sbin/revnetgroup
sbin/rpc.yppasswdd
sbin/rpc.ypxfrd
sbin/yphelper
sbin/ypinit
sbin/yppush
sbin/ypserv
sbin/ypxfr
sbin/ypxfr_1perday
sbin/ypxfr_1perhour
sbin/ypxfr_2perday
share/examples/rc.d/yppasswdd
share/examples/rc.d/ypserv
share/man/man5/netgroup.5
share/man/man5/ypserv.conf.5
share/man/man8/makedbm.8
share/man/man8/mknetid.8
share/man/man8/pwupdate.8
share/man/man8/revnetgroup.8
share/man/man8/rpc.yppasswdd.8
share/man/man8/rpc.ypxfrd.8
share/man/man8/yphelper.8
share/man/man8/ypinit.8
share/man/man8/yppasswdd.8
share/man/man8/yppush.8
share/man/man8/ypserv.8
share/man/man8/ypxfr.8
share/man/man8/ypxfrd.8
