# thy/PLIST
@comment $NetBSD: PLIST,v 1.4 2004/03/21 14:19:02 algernon Exp $
etc/rc.d/thy
share/thy/samples/http301.html
share/thy/samples/http302.html
share/thy/samples/http400.html
share/thy/samples/http401.html
share/thy/samples/http403.html
share/thy/samples/http404.html
share/thy/samples/http405.html
share/thy/samples/http411.html
share/thy/samples/http412.html
share/thy/samples/http413.html
share/thy/samples/http416.html
share/thy/samples/http417.html
share/thy/samples/http500.html
share/thy/samples/http501.html
share/thy/samples/http503.html
share/thy/samples/http505.html
share/thy/samples/thyrc
@dirrm share/thy/samples
@dirrm share/thy
info/thy.info
man/man8/thy.8
man/man7/thy-auth.7
sbin/thy
