# drweb/files/drwebd-rcd.src
#!@SH@
#
# $NetBSD: drwebd-rcd.src,v 1.1 2003/09/11 15:20:28 bubuchka Exp $

# PROVIDE: drwebd
# REQUIRE: SERVERS
# BEFORE: DAEMON

. /etc/rc.subr

name="drwebd"
rcvar=$name
command="@PREFIX@/sbin/${name}-bin"
procname="${command}"
drw_inifile="@PKG_SYSCONFDIR@/drweb32.ini"
required_files="${drw_inifile}"
command_args="-ini:${drw_inifile}"

extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
