# mu-conference/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/10/11 10:44:58 thomasklausner Exp $
bin/mu-conference
share/doc/muc/ChangeLog
share/doc/muc/FAQ
share/doc/muc/README
share/doc/muc/TODO
share/doc/muc/XEP0045_SUPPORT
share/examples/jabber/muc-default.xml
share/examples/jabber/style.css
share/examples/rc.d/muc
