# rsyslog5/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2011/04/18 13:46:22 bartoszkuzma Exp $
${PLIST.gnutls}lib/rsyslog/lmnsd_gtls.la
${PLIST.gssapi}lib/rsyslog/imgssapi.la
${PLIST.gssapi}lib/rsyslog/lmgssutil.la
${PLIST.gssapi}lib/rsyslog/omgssapi.la
${PLIST.mail}lib/rsyslog/ommail.la
${PLIST.mysql}lib/rsyslog/ommysql.la
${PLIST.mysql}share/examples/rsyslog/createDB-mysql.sql
${PLIST.pgsql}lib/rsyslog/ompgsql.la
${PLIST.pgsql}share/examples/rsyslog/createDB-pgsql.sql
${PLIST.snmp}lib/rsyslog/omsnmp.la
lib/rsyslog/imklog.la
lib/rsyslog/immark.la
lib/rsyslog/imtcp.la
lib/rsyslog/imudp.la
lib/rsyslog/imuxsock.la
lib/rsyslog/lmnet.la
lib/rsyslog/lmnetstrms.la
lib/rsyslog/lmnsd_ptcp.la
lib/rsyslog/lmregexp.la
lib/rsyslog/lmstrmsrv.la
lib/rsyslog/lmtcpclt.la
lib/rsyslog/lmtcpsrv.la
lib/rsyslog/lmzlibw.la
lib/rsyslog/omruleset.la
lib/rsyslog/omtesting.la
man/man5/rsyslog.conf.5
man/man8/rsyslogd.8
sbin/rsyslogd
share/examples/rc.d/rsyslogd
share/examples/rsyslog/rsyslog.conf
