# rsyslog5/PLIST
@comment $NetBSD: PLIST,v 1.2 2011/10/27 10:14:19 fhajny Exp $
${PLIST.file}lib/rsyslog/imfile.la
${PLIST.gssapi}lib/rsyslog/imgssapi.la
${PLIST.klog}lib/rsyslog/imklog.la
lib/rsyslog/immark.la
${PLIST.solaris}lib/rsyslog/imsolaris.la
lib/rsyslog/imtcp.la
lib/rsyslog/imudp.la
lib/rsyslog/imuxsock.la
${PLIST.gssapi}lib/rsyslog/lmgssutil.la
lib/rsyslog/lmnet.la
lib/rsyslog/lmnetstrms.la
${PLIST.gnutls}lib/rsyslog/lmnsd_gtls.la
lib/rsyslog/lmnsd_ptcp.la
lib/rsyslog/lmregexp.la
lib/rsyslog/lmstrmsrv.la
lib/rsyslog/lmtcpclt.la
lib/rsyslog/lmtcpsrv.la
lib/rsyslog/lmzlibw.la
${PLIST.gssapi}lib/rsyslog/omgssapi.la
${PLIST.mail}lib/rsyslog/ommail.la
${PLIST.mysql}lib/rsyslog/ommysql.la
${PLIST.pgsql}lib/rsyslog/ompgsql.la
lib/rsyslog/omruleset.la
${PLIST.snmp}lib/rsyslog/omsnmp.la
lib/rsyslog/omtesting.la
man/man5/rsyslog.conf.5
man/man8/rsyslogd.8
sbin/rsyslogd
share/examples/rc.d/rsyslogd
${PLIST.mysql}share/examples/rsyslog/createDB-mysql.sql
${PLIST.pgsql}share/examples/rsyslog/createDB-pgsql.sql
share/examples/rsyslog/rsyslog.conf
