# adguardhome/PLIST
@comment $NetBSD$
sbin/adguardhome
share/examples/rc.d/adguardhome
@pkgdir libdata/adguardhome
@pkgdir etc/adguardhome
