# quickml/PLIST
@comment $NetBSD: PLIST,v 1.2 2011/05/08 01:58:04 obache Exp $
bin/quickml-analog
${RUBY_SITELIB}/quickml.rb
${RUBY_SITELIB}/quickml/config.rb
${RUBY_SITELIB}/quickml/core.rb
${RUBY_SITELIB}/quickml/gettext.rb
${RUBY_SITELIB}/quickml/logger.rb
${RUBY_SITELIB}/quickml/mail.rb
${RUBY_SITELIB}/quickml/server.rb
${RUBY_SITELIB}/quickml/sweeper.rb
${RUBY_SITELIB}/quickml/utils.rb
${RUBY_SITELIB}/quickml/version.rb
sbin/quickml
sbin/quickml-ctl
share/examples/quickml/quickmlrc.sample
share/examples/rc.d/quickml
