# linuxigd/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:00 jeremy-c-reed Exp $
bin/upnpd
share/examples/linuxigd/gateconnSCPD.xml
share/examples/linuxigd/gatedesc.skl
share/examples/linuxigd/gatedesc.xml
share/examples/linuxigd/gateicfgSCPD.xml
share/examples/linuxigd/gateinfoSCPD.xml
share/examples/rc.d/linuxigd
@dirrm share/examples/linuxigd
