# linuxigd/PLIST
@comment $NetBSD: PLIST,v 1.3 2006/08/27 15:00:00 obache Exp $
man/man8/upnpd.8
sbin/upnpd
etc/linuxigd/dummy.xml
etc/linuxigd/gateconnSCPD.xml
etc/linuxigd/gatedesc.xml
etc/linuxigd/gateicfgSCPD.xml
etc/linuxigd/ligd.gif
share/examples/upnpd.conf
share/examples/rc.d/linuxigd
