# munge/PLIST.Linux
@comment $NetBSD: PLIST.Linux,v 1.4 2018/03/24 18:05:58 bacon Exp $
etc/logrotate.d/munge
etc/sysconfig/munge
lib/systemd/system/munge.service
