# munge/PLIST.Linux
@comment $NetBSD: PLIST.Linux,v 1.4 2014/04/29 19:38:51 outpaddling Exp $
share/examples/rc.d/init.d/munge
share/examples/sysconfig/munge
