# linuxigd-current/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2005/10/02 13:35:54 obache Exp $
sbin/upnpd
share/examples/linuxigd/upnpd.conf
share/examples/linuxigd/dummy.xml
share/examples/linuxigd/gateconnSCPD.xml
share/examples/linuxigd/gatedesc.xml
share/examples/linuxigd/gateicfgSCPD.xml
share/examples/linuxigd/ligd.gif
share/examples/rc.d/linuxigd
@dirrm share/examples/linuxigd
