# dkim-milter/PLIST
@comment $NetBSD: PLIST,v 1.4 2007/08/02 23:33:43 jukka Exp $
lib/libdkim.a
man/man5/dkim-milter.conf.5
man/man8/dkim-milter.8
sbin/dkim-milter
share/doc/dkim-milter/FEATURES
share/doc/dkim-milter/INSTALL
share/doc/dkim-milter/KNOWNBUGS
share/doc/dkim-milter/LICENSE
share/doc/dkim-milter/README
share/doc/dkim-milter/RELEASE_NOTES
share/doc/dkim-milter/draft-hoffman-dac-vbr-01.txt
share/doc/dkim-milter/draft-ietf-dkim-overview-05.txt
share/doc/dkim-milter/draft-ietf-dkim-ssp-00.txt
share/doc/dkim-milter/draft-ietf-dkim-ssp-requirements-04.txt
share/doc/dkim-milter/rfc4686.txt
share/doc/dkim-milter/rfc4871.txt
share/examples/dkim-milter/dkim-milter.conf.sample
share/examples/rc.d/dkim-milter
@dirrm share/doc/dkim-milter
@dirrm share/examples/dkim-milter
