# libreswan/PLIST
@comment $NetBSD$
libexec/ipsec/_import_crl
libexec/ipsec/_plutorun
libexec/ipsec/_realsetup
libexec/ipsec/_secretcensor
libexec/ipsec/_unbound-hook
libexec/ipsec/_updown
libexec/ipsec/_updown.bsd
libexec/ipsec/addconn
libexec/ipsec/algparse
libexec/ipsec/asn1check
libexec/ipsec/auto
libexec/ipsec/barf
libexec/ipsec/cavp
libexec/ipsec/dncheck
libexec/ipsec/ecdsasigkey
libexec/ipsec/enumcheck
libexec/ipsec/hunkcheck
libexec/ipsec/ipcheck
libexec/ipsec/jambufcheck
libexec/ipsec/keyidcheck
libexec/ipsec/letsencrypt
libexec/ipsec/newhostkey
libexec/ipsec/pluto
libexec/ipsec/readwriteconf
libexec/ipsec/rsasigkey
libexec/ipsec/setup
libexec/ipsec/showhostkey
libexec/ipsec/showroute
libexec/ipsec/timecheck
libexec/ipsec/vendoridcheck
libexec/ipsec/whack
man/man5/ipsec.conf.5
man/man5/ipsec.secrets.5
man/man8/ipsec.8
man/man8/ipsec__import_crl.8
man/man8/ipsec__plutorun.8
man/man8/ipsec__realsetup.8
man/man8/ipsec__secretcensor.8
man/man8/ipsec__unbound-hook.8
man/man8/ipsec__updown.8
man/man8/ipsec__updown.bsdkame.8
man/man8/ipsec_addconn.8
man/man8/ipsec_auto.8
man/man8/ipsec_barf.8
man/man8/ipsec_checknss.8
man/man8/ipsec_ecdsasigkey.8
man/man8/ipsec_import.8
man/man8/ipsec_initnss.8
man/man8/ipsec_letsencrypt.8
man/man8/ipsec_newhostkey.8
man/man8/ipsec_pluto.8
man/man8/ipsec_readwriteconf.8
man/man8/ipsec_rsasigkey.8
man/man8/ipsec_setup.8
man/man8/ipsec_showhostkey.8
man/man8/ipsec_showroute.8
man/man8/ipsec_vendorid.8
man/man8/ipsec_whack.8
man/man8/pluto.8
sbin/ipsec
share/examples/libreswan/ipsec.conf-sample
share/examples/libreswan/ipsec.d/policies/block
share/examples/libreswan/ipsec.d/policies/clear
share/examples/libreswan/ipsec.d/policies/clear-or-private
share/examples/libreswan/ipsec.d/policies/portexcludes.conf
share/examples/libreswan/ipsec.d/policies/private
share/examples/libreswan/ipsec.d/policies/private-or-clear
share/examples/libreswan/ipsec.secrets-sample
share/examples/libreswan/logrotate.d/libreswan
share/examples/libreswan/pam.d/pluto
share/examples/libreswan/rc.d/pluto.sh
