# isc-dhclient/PLIST
@comment $NetBSD: PLIST,v 1.2 2004/04/18 17:22:12 mpasternak Exp $
etc/rc.d/isc_dhclient
man/man5/dhclient.conf.5
man/man5/dhclient.leases.5
man/man8/dhclient-script.8
man/man8/dhclient.8
sbin/dhclient
sbin/dhclient-script
share/doc/isc-dhclient/README
@dirrm share/doc/isc-dhclient
