# isc-dhclient/PLIST
@comment $NetBSD: PLIST,v 1.6 2007/12/13 06:47:43 marttikuparinen Exp $
man/man5/dhclient.conf.5
man/man5/dhclient.leases.5
man/man8/dhclient-script.8
man/man8/dhclient.8
sbin/dhclient
sbin/dhclient-script
share/examples/isc-dhcp-client/dhclient.conf
share/examples/rc.d/isc_dhclient
@dirrm share/examples/isc-dhcp-client
