# cvsd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2003/04/05 15:21:05 virtus Exp $
etc/init.d/cvsd
man/man5/cvsd.conf.5
man/man8/cvsd-buildroot.8
man/man8/cvsd-passwd.8
man/man8/cvsd.8
sbin/cvsd
sbin/cvsd-buildroot
sbin/cvsd-passwd
@dirrm etc/init.d
