# cvsd/PLIST
@comment $NetBSD: PLIST,v 1.5 2004/07/28 21:50:53 schnoebe Exp $
etc/rc.d/cvsd
sbin/cvsd
sbin/cvsd-buildroot
sbin/cvsd-passwd
man/man5/cvsd.conf.5
man/man8/cvsd-buildroot.8
man/man8/cvsd-passwd.8
man/man8/cvsd.8
share/doc/cvsd/README
share/examples/cvsd/cvsd.conf.example
@dirrm share/examples/cvsd
@dirrm share/doc/cvsd
