# cvsd/PLIST
@comment $NetBSD: PLIST,v 1.2 2003/04/20 16:14:35 virtus Exp $
etc/rc.d/cvsd
man/man5/cvsd.conf.5
man/man8/cvsd-buildroot.8
man/man8/cvsd-passwd.8
man/man8/cvsd.8
sbin/cvsd
sbin/cvsd-buildroot
sbin/cvsd-passwd
