# netkit-base/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:01 jeremy-c-reed Exp $
bin/ping
man/man8/inetd.8
man/man8/ping.8
sbin/inetd
share/examples/${PKGBASE}/etc.sample/README
share/examples/${PKGBASE}/etc.sample/diphosts
share/examples/${PKGBASE}/etc.sample/exports
share/examples/${PKGBASE}/etc.sample/ftpusers
share/examples/${PKGBASE}/etc.sample/host.conf
share/examples/${PKGBASE}/etc.sample/hosts
share/examples/${PKGBASE}/etc.sample/hosts.allow
share/examples/${PKGBASE}/etc.sample/hosts.deny
share/examples/${PKGBASE}/etc.sample/hosts.equiv
share/examples/${PKGBASE}/etc.sample/hosts.lpd
share/examples/${PKGBASE}/etc.sample/inetd.conf
share/examples/${PKGBASE}/etc.sample/networks
share/examples/${PKGBASE}/etc.sample/nntpserver
share/examples/${PKGBASE}/etc.sample/nsswitch.conf
share/examples/${PKGBASE}/etc.sample/printcap
share/examples/${PKGBASE}/etc.sample/protocols
share/examples/${PKGBASE}/etc.sample/resolv.conf
share/examples/${PKGBASE}/etc.sample/rpc
share/examples/${PKGBASE}/etc.sample/sample.dip
share/examples/${PKGBASE}/etc.sample/services
share/examples/${PKGBASE}/etc.sample/slip.hosts
share/examples/${PKGBASE}/etc.sample/slip.login
share/examples/${PKGBASE}/inetd.conf
share/examples/rc.d/inetd
@dirrm share/examples/${PKGBASE}/etc.sample
@dirrm share/examples/${PKGBASE}
