# stud/PLIST
@comment $NetBSD: PLIST,v 1.3 2013/03/07 14:21:53 jym-netbsd Exp $
bin/stud
man/man8/stud.8
share/examples/stud/stud.conf
share/examples/rc.d/stud
