# ja-freewnn-cvs/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2011/04/28 13:48:14 makoto Exp $
bin/Wnn4/atod
bin/Wnn4/atof
bin/Wnn4/dtoa
bin/Wnn4/jserver
bin/Wnn4/oldatonewa
bin/Wnn4/wddel
bin/Wnn4/wdreg
bin/Wnn4/wnnkill
bin/Wnn4/wnnstat
bin/Wnn4/wnntouch
bin/cWnn4/catod
bin/cWnn4/catof
bin/cWnn4/cdtoa
bin/cWnn4/cserver
bin/cWnn4/cwddel
bin/cWnn4/cwdreg
bin/cWnn4/cwnnkill
bin/cWnn4/cwnnstat
bin/cWnn4/cwnntouch
bin/cWnn4/tserver
bin/kWnn4/katod
bin/kWnn4/katof
bin/kWnn4/kdtoa
bin/kWnn4/kserver
bin/kWnn4/kwddel
bin/kWnn4/kwdreg
bin/kWnn4/kwnnkill
bin/kWnn4/kwnnstat
bin/kWnn4/kwnntouch
include/cwnn/commonhd.h
include/cwnn/cplib.h
include/cwnn/jllib.h
include/cwnn/jslib.h
include/cwnn/msg.h
include/cwnn/rk_spclval.h
include/cwnn/wnnerror.h
include/kwnn/commonhd.h
include/kwnn/cplib.h
include/kwnn/jllib.h
include/kwnn/jslib.h
include/kwnn/msg.h
include/kwnn/rk_spclval.h
include/kwnn/wnnerror.h
include/wnn/commonhd.h
include/wnn/config.h
include/wnn/cplib.h
include/wnn/jlib.h
include/wnn/jllib.h
include/wnn/jslib.h
include/wnn/msg.h
include/wnn/rk_spclval.h
include/wnn/wnn_config.h
include/wnn/wnnerror.h
lib/libcwnn.la
lib/libjd.la
lib/libkwnn.la
lib/libwnn.la
lib/wnn/cvt_key_empty
lib/wnn/cvt_key_tbl
lib/wnn/cvt_key_tbl.ST
lib/wnn/cvt_key_tbl.gm
lib/wnn/cvt_key_tbl.kt
lib/wnn/cvt_key_tbl.mv
lib/wnn/cvt_key_tbl.vt
lib/wnn/ja_JP/dic/gerodic/g-jinmei.dic
lib/wnn/ja_JP/dic/pubdic/bio.dic
lib/wnn/ja_JP/dic/pubdic/chimei.dic
lib/wnn/ja_JP/dic/pubdic/computer.dic
lib/wnn/ja_JP/dic/pubdic/full.fzk
lib/wnn/ja_JP/dic/pubdic/jinmei.dic
lib/wnn/ja_JP/dic/pubdic/kihon.dic
lib/wnn/ja_JP/dic/pubdic/kougo.fzk
lib/wnn/ja_JP/dic/pubdic/koyuu.dic
lib/wnn/ja_JP/dic/pubdic/setsuji.dic
lib/wnn/ja_JP/dic/pubdic/special.dic
lib/wnn/ja_JP/dic/pubdic/std.fzk
lib/wnn/ja_JP/dic/pubdic/symbol.dic
lib/wnn/ja_JP/dic/pubdic/tankan.dic
lib/wnn/ja_JP/dic/src/fzk.attr
lib/wnn/ja_JP/dic/src/fzk.con
lib/wnn/ja_JP/dic/src/fzk.fzkattr
lib/wnn/ja_JP/dic/src/fzk.jirattr
lib/wnn/ja_JP/dic/src/fzk.jircon
lib/wnn/ja_JP/dic/src/fzk.master
lib/wnn/ja_JP/dic/src/fzk.shuutan
lib/wnn/ja_JP/dic/wnncons/tankan2.dic
lib/wnn/ja_JP/dic/wnncons/tankan3.dic
lib/wnn/ja_JP/hinsi.data
lib/wnn/ja_JP/jserverrc
lib/wnn/ja_JP/libwnn.msg
lib/wnn/ja_JP/rk.vi/1B_newTOUPPER
lib/wnn/ja_JP/rk.vi/2A_CTRL
lib/wnn/ja_JP/rk.vi/2B_KEISEN
lib/wnn/ja_JP/rk.vi/2B_ROMKANA
lib/wnn/ja_JP/rk.vi/2C_KEISEN1
lib/wnn/ja_JP/rk.vi/2C_VI
lib/wnn/ja_JP/rk.vi/2C_VIEX
lib/wnn/ja_JP/rk.vi/2C_VISH
lib/wnn/ja_JP/rk.vi/2_VITHROW
lib/wnn/ja_JP/rk.vi/3B_ZENASC
lib/wnn/ja_JP/rk.vi/mode
lib/wnn/ja_JP/rk.vi/uumkey
lib/wnn/ja_JP/rk/1B_TOUPPER
lib/wnn/ja_JP/rk/1B_ZENHIRA
lib/wnn/ja_JP/rk/1B_ZENKATA
lib/wnn/ja_JP/rk/1B_newTOUPPER
lib/wnn/ja_JP/rk/2A_CTRL
lib/wnn/ja_JP/rk/2B_DAKUTEN
lib/wnn/ja_JP/rk/2B_JIS
lib/wnn/ja_JP/rk/2B_ROMKANA
lib/wnn/ja_JP/rk/2_TCODE
lib/wnn/ja_JP/rk/3B_HANKATA
lib/wnn/ja_JP/rk/3B_KATAKANA
lib/wnn/ja_JP/rk/3B_ZENKAKU
lib/wnn/ja_JP/rk/autork
lib/wnn/ja_JP/rk/mode
lib/wnn/ja_JP/rk/mode.hankata
lib/wnn/ja_JP/rk/mode.nohankata
lib/wnn/ja_JP/uum.msg
lib/wnn/ja_JP/uumkey
lib/wnn/ja_JP/uumkey.omr
lib/wnn/ja_JP/uumkey_e
lib/wnn/ja_JP/uumrc
lib/wnn/ja_JP/uumrc.omr
lib/wnn/ja_JP/uumrc.rev
lib/wnn/ja_JP/uumrc_e
lib/wnn/ja_JP/uumrc_vi
lib/wnn/ja_JP/wnnenvrc
lib/wnn/ja_JP/wnnenvrc.omr
lib/wnn/ja_JP/wnnenvrc.rem
lib/wnn/ja_JP/wnnenvrc.rev
lib/wnn/ja_JP/wnnenvrc_R
lib/wnn/ja_JP/wnnenvrc_R.omr
lib/wnn/ja_JP/wnnenvrc_R.rem
lib/wnn/ja_JP/wnnstat.msg
lib/wnn/ko_KR/dic/sys/full.fzk
lib/wnn/ko_KR/dic/sys/hword.dic
lib/wnn/ko_KR/dic/sys/single.dic
lib/wnn/ko_KR/hinsi.data
lib/wnn/ko_KR/kserverrc
lib/wnn/ko_KR/libwnn.msg
lib/wnn/ko_KR/rk/1B_TOUPPER
lib/wnn/ko_KR/rk/2A_CTRL
lib/wnn/ko_KR/rk/2C_HANGUL.KEY
lib/wnn/ko_KR/rk/2C_HANGUL.ROM
lib/wnn/ko_KR/rk/2C_KSC
lib/wnn/ko_KR/rk/3B_JEONKAK
lib/wnn/ko_KR/rk/mode
lib/wnn/ko_KR/rk/mode.key
lib/wnn/ko_KR/rk/mode.rom
lib/wnn/ko_KR/uum.msg
lib/wnn/ko_KR/uumkey
lib/wnn/ko_KR/uumrc
lib/wnn/ko_KR/wnnenvrc
lib/wnn/ko_KR/wnnenvrc_R
lib/wnn/ko_KR/wnnstat.msg
lib/wnn/lt_LN/rk/2A_CTRL
lib/wnn/lt_LN/rk/2B_LATIN
lib/wnn/lt_LN/rk/mode
lib/wnn/lt_LN/uum.msg
lib/wnn/lt_LN/uumkey
lib/wnn/lt_LN/uumkey_e
lib/wnn/lt_LN/uumrc
lib/wnn/serverdefs
lib/wnn/zh_CN/cixing.data
lib/wnn/zh_CN/cserverrc
lib/wnn/zh_CN/dic/sys/QianMa.dic
lib/wnn/zh_CN/dic/sys/WuBi.dic
lib/wnn/zh_CN/dic/sys/basic.dic
lib/wnn/zh_CN/dic/sys/computer.dic
lib/wnn/zh_CN/dic/sys/cwnn.dic
lib/wnn/zh_CN/dic/sys/full.con
lib/wnn/zh_CN/dic/sys/full.conR
lib/wnn/zh_CN/dic/sys/level_1.dic
lib/wnn/zh_CN/dic/sys/level_2.dic
lib/wnn/zh_CN/libwnn.msg
lib/wnn/zh_CN/rk/1B_BS
lib/wnn/zh_CN/rk/1B_TOLOWER
lib/wnn/zh_CN/rk/1B_TOUPPER
lib/wnn/zh_CN/rk/1Z_ZhuYin
lib/wnn/zh_CN/rk/2A_CTRL
lib/wnn/zh_CN/rk/2B_WuBi
lib/wnn/zh_CN/rk/2P_QuanPin
lib/wnn/zh_CN/rk/2P_RongCuo
lib/wnn/zh_CN/rk/2Z_ZhuYin
lib/wnn/zh_CN/rk/2Z_tail
lib/wnn/zh_CN/rk/2Z_tail_ma
lib/wnn/zh_CN/rk/2Z_tail_pin
lib/wnn/zh_CN/rk/3B_quanjiao
lib/wnn/zh_CN/rk/mode
lib/wnn/zh_CN/rk_p/1B_BS
lib/wnn/zh_CN/rk_p/1B_TOLOWER
lib/wnn/zh_CN/rk_p/1B_TOUPPER
lib/wnn/zh_CN/rk_p/1F_select
lib/wnn/zh_CN/rk_p/2A_CTRL
lib/wnn/zh_CN/rk_p/2B_QianMa
lib/wnn/zh_CN/rk_p/2B_WuBi
lib/wnn/zh_CN/rk_p/2B_ch_yo_zi
lib/wnn/zh_CN/rk_p/2B_ma_GuoBi
lib/wnn/zh_CN/rk_p/2B_ma_NeiMa
lib/wnn/zh_CN/rk_p/2B_ma_QuWei
lib/wnn/zh_CN/rk_p/2F_restart
lib/wnn/zh_CN/rk_p/2F_select
lib/wnn/zh_CN/rk_p/2P_ErPin
lib/wnn/zh_CN/rk_p/2P_ErPin_3K
lib/wnn/zh_CN/rk_p/2P_ErPin_P
lib/wnn/zh_CN/rk_p/2P_ErPin_P3K
lib/wnn/zh_CN/rk_p/2P_QuanPin
lib/wnn/zh_CN/rk_p/2P_RC2
lib/wnn/zh_CN/rk_p/2P_RC2_P
lib/wnn/zh_CN/rk_p/2P_RongCuo
lib/wnn/zh_CN/rk_p/2Z_tail_ma
lib/wnn/zh_CN/rk_p/2Z_tail_pin
lib/wnn/zh_CN/rk_p/3B_quanjiao
lib/wnn/zh_CN/rk_p/3P_ErPin
lib/wnn/zh_CN/rk_p/3P_ErPin_P
lib/wnn/zh_CN/rk_p/mode
lib/wnn/zh_CN/rk_z/1B_TOUPPER
lib/wnn/zh_CN/rk_z/1F_select
lib/wnn/zh_CN/rk_z/1Z_ZhuYin
lib/wnn/zh_CN/rk_z/1Z_ZhuYin_X
lib/wnn/zh_CN/rk_z/2A_CTRL
lib/wnn/zh_CN/rk_z/2F_select
lib/wnn/zh_CN/rk_z/2Z_ZhuYin
lib/wnn/zh_CN/rk_z/2Z_ZhuYin_X
lib/wnn/zh_CN/rk_z/2Z_tail
lib/wnn/zh_CN/rk_z/2Z_tail_X
lib/wnn/zh_CN/rk_z/3B_quanjiao
lib/wnn/zh_CN/rk_z/mode
lib/wnn/zh_CN/uum.msg
lib/wnn/zh_CN/uumkey
lib/wnn/zh_CN/uumkey_e
lib/wnn/zh_CN/uumkey_p
lib/wnn/zh_CN/uumrc
lib/wnn/zh_CN/uumrc_p
lib/wnn/zh_CN/uumrc_z
lib/wnn/zh_CN/wnnenvrc
lib/wnn/zh_CN/wnnenvrc_Qi
lib/wnn/zh_CN/wnnenvrc_QiR
lib/wnn/zh_CN/wnnenvrc_R
lib/wnn/zh_CN/wnnenvrc_Wu
lib/wnn/zh_CN/wnnenvrc_WuR
lib/wnn/zh_CN/wnnstat.msg
lib/wnn/zh_TW/cixing.data
lib/wnn/zh_TW/dic/sys/cns_ch.dic
lib/wnn/zh_TW/dic/sys/cns_wd.dic
lib/wnn/zh_TW/dic/sys/full.con
lib/wnn/zh_TW/dic/sys/full.conR
lib/wnn/zh_TW/libwnn.msg
lib/wnn/zh_TW/rk/2A_CTRL
lib/wnn/zh_TW/rk/mode
lib/wnn/zh_TW/rk_p/2A_CTRL
lib/wnn/zh_TW/rk_p/mode
lib/wnn/zh_TW/rk_z/2A_CTRL
lib/wnn/zh_TW/rk_z/mode
lib/wnn/zh_TW/tserverrc
lib/wnn/zh_TW/uum.msg
lib/wnn/zh_TW/uumkey
lib/wnn/zh_TW/uumkey_e
lib/wnn/zh_TW/uumkey_p
lib/wnn/zh_TW/uumrc
lib/wnn/zh_TW/uumrc_p
lib/wnn/zh_TW/uumrc_z
lib/wnn/zh_TW/wnnenvrc
lib/wnn/zh_TW/wnnenvrc_R
lib/wnn/zh_TW/wnnstat.msg
share/examples/rc.d/ja_freewnn
