# sysklogd/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:34:06 jeremy-c-reed Exp $
man/man5/syslog.conf.5
man/man8/klogd.8
man/man8/sysklogd.8
man/man8/syslogd.8
sbin/klogd
sbin/syslogd
share/examples/rc.d/klogd
share/examples/rc.d/syslogd
share/examples/sysklogd/syslog.conf
share/examples/sysklogd/syslog.conf.example
@dirrm share/examples/sysklogd
