# 389-ds-base/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/09/01 21:21:31 udontknow Exp $
bin/cl-dump
bin/cl-dump.pl
bin/dbgen.pl
bin/dbscan
bin/dbscan-bin
bin/dsktune
bin/dsktune-bin
bin/infadd
bin/infadd-bin
bin/ldap-agent
bin/ldap-agent-bin
bin/ldclt
bin/ldclt-bin
bin/ldif
bin/ldif-bin
bin/logconv.pl
bin/migratecred
bin/migratecred-bin
bin/mmldif
bin/mmldif-bin
bin/pwdhash
bin/pwdhash-bin
bin/repl-monitor
bin/repl-monitor.pl
bin/rsearch
bin/rsearch-bin
lib/dirsrv/libns-dshttpd.la
lib/dirsrv/libslapd.la
lib/dirsrv/perl/DSCreate.pm
lib/dirsrv/perl/DSDialogs.pm
lib/dirsrv/perl/DSMigration.pm
lib/dirsrv/perl/Dialog.pm
lib/dirsrv/perl/DialogManager.pm
lib/dirsrv/perl/FileConn.pm
lib/dirsrv/perl/Inf.pm
lib/dirsrv/perl/Migration.pm
lib/dirsrv/perl/Resource.pm
lib/dirsrv/perl/Setup.pm
lib/dirsrv/perl/SetupDialogs.pm
lib/dirsrv/perl/SetupLog.pm
lib/dirsrv/perl/Util.pm
lib/dirsrv/plugins/libacl-plugin.la
lib/dirsrv/plugins/libattr-unique-plugin.la
lib/dirsrv/plugins/libback-ldbm.la
lib/dirsrv/plugins/libbitwise-plugin.la
lib/dirsrv/plugins/libchainingdb-plugin.la
lib/dirsrv/plugins/libcollation-plugin.la
lib/dirsrv/plugins/libcos-plugin.la
lib/dirsrv/plugins/libderef-plugin.la
lib/dirsrv/plugins/libdes-plugin.la
lib/dirsrv/plugins/libdistrib-plugin.la
lib/dirsrv/plugins/libdna-plugin.la
lib/dirsrv/plugins/libhttp-client-plugin.la
lib/dirsrv/plugins/liblinkedattrs-plugin.la
lib/dirsrv/plugins/libmemberof-plugin.la
lib/dirsrv/plugins/libpam-passthru-plugin.la
lib/dirsrv/plugins/libpassthru-plugin.la
lib/dirsrv/plugins/libpwdstorage-plugin.la
lib/dirsrv/plugins/libreferint-plugin.la
lib/dirsrv/plugins/libreplication-plugin.la
lib/dirsrv/plugins/libretrocl-plugin.la
lib/dirsrv/plugins/libroles-plugin.la
lib/dirsrv/plugins/libschemareload-plugin.la
lib/dirsrv/plugins/libstatechange-plugin.la
lib/dirsrv/plugins/libsyntax-plugin.la
lib/dirsrv/plugins/libusn-plugin.la
lib/dirsrv/plugins/libviews-plugin.la
man/man1/cl-dump.1
man/man1/dbgen.pl.1
man/man1/dbscan.1
man/man1/dsktune.1
man/man1/infadd.1
man/man1/ldap-agent.1
man/man1/ldclt.1
man/man1/ldif.1
man/man1/logconv.pl.1
man/man1/migratecred.1
man/man1/mmldif.1
man/man1/pwdhash.1
man/man1/repl-monitor.1
man/man1/rsearch.1
man/man8/migrate-ds.pl.8
man/man8/ns-slapd.8
man/man8/remove-ds.pl.8
man/man8/setup-ds.pl.8
sbin/migrate-ds.pl
sbin/ns-slapd
sbin/remove-ds.pl
sbin/setup-ds.pl
share/dirsrv/data/10rfc2307bis.ldif
share/dirsrv/data/60changelog.ldif
share/dirsrv/data/60inetmail.ldif
share/dirsrv/data/60kerberos.ldif
share/dirsrv/data/60krb5kdc.ldif
share/dirsrv/data/60nis.ldif
share/dirsrv/data/60qmail.ldif
share/dirsrv/data/60radius.ldif
share/dirsrv/data/60rfc4876.ldif
share/dirsrv/data/60samba.ldif
share/dirsrv/data/60samba3.ldif
share/dirsrv/data/Ace.ldif
share/dirsrv/data/European.ldif
share/dirsrv/data/Eurosuffix.ldif
share/dirsrv/data/Example-roles.ldif
share/dirsrv/data/Example-views.ldif
share/dirsrv/data/Example.ldif
share/dirsrv/data/dbgen-FamilyNames
share/dirsrv/data/dbgen-GivenNames
share/dirsrv/data/dbgen-OrgUnits
share/dirsrv/data/template-baseacis.ldif
share/dirsrv/data/template-bitwise.ldif
share/dirsrv/data/template-country.ldif
share/dirsrv/data/template-dnaplugin.ldif
share/dirsrv/data/template-domain.ldif
share/dirsrv/data/template-dse.ldif
share/dirsrv/data/template-ldapi-autobind.ldif
share/dirsrv/data/template-ldapi-default.ldif
share/dirsrv/data/template-ldapi.ldif
share/dirsrv/data/template-locality.ldif
share/dirsrv/data/template-org.ldif
share/dirsrv/data/template-orgunit.ldif
share/dirsrv/data/template-pampta.ldif
share/dirsrv/data/template-sasl.ldif
share/dirsrv/data/template-state.ldif
share/dirsrv/data/template-suffix-db.ldif
share/dirsrv/data/template.ldif
share/dirsrv/inf/dscreate.map
share/dirsrv/inf/dsorgentries.map
share/dirsrv/inf/slapd.inf
share/dirsrv/mibs/NETWORK-SERVICES-MIB.txt
share/dirsrv/mibs/RFC-1215.txt
share/dirsrv/mibs/RFC1155-SMI.txt
share/dirsrv/mibs/SNMPv2-CONF.txt
share/dirsrv/mibs/SNMPv2-SMI.txt
share/dirsrv/mibs/SNMPv2-TC.txt
share/dirsrv/mibs/redhat-directory.mib
share/dirsrv/properties/migrate-ds.res
share/dirsrv/properties/ns-slapd.properties
share/dirsrv/properties/setup-ds.res
share/dirsrv/script-templates/template-bak2db
share/dirsrv/script-templates/template-bak2db.pl
share/dirsrv/script-templates/template-db2bak
share/dirsrv/script-templates/template-db2bak.pl
share/dirsrv/script-templates/template-db2index
share/dirsrv/script-templates/template-db2index.pl
share/dirsrv/script-templates/template-db2ldif
share/dirsrv/script-templates/template-db2ldif.pl
share/dirsrv/script-templates/template-dbverify
share/dirsrv/script-templates/template-fixup-linkedattrs.pl
share/dirsrv/script-templates/template-fixup-memberof.pl
share/dirsrv/script-templates/template-ldif2db
share/dirsrv/script-templates/template-ldif2db.pl
share/dirsrv/script-templates/template-ldif2ldap
share/dirsrv/script-templates/template-monitor
share/dirsrv/script-templates/template-ns-accountstatus.pl
share/dirsrv/script-templates/template-ns-activate.pl
share/dirsrv/script-templates/template-ns-inactivate.pl
share/dirsrv/script-templates/template-ns-newpwpolicy.pl
share/dirsrv/script-templates/template-restart-slapd
share/dirsrv/script-templates/template-restoreconfig
share/dirsrv/script-templates/template-saveconfig
share/dirsrv/script-templates/template-schema-reload.pl
share/dirsrv/script-templates/template-start-slapd
share/dirsrv/script-templates/template-stop-slapd
share/dirsrv/script-templates/template-suffix2instance
share/dirsrv/script-templates/template-syntax-validate.pl
share/dirsrv/script-templates/template-usn-tombstone-cleanup.pl
share/dirsrv/script-templates/template-verify-db.pl
share/dirsrv/script-templates/template-vlvindex
share/examples/389-ds-base/dirsrv/config/certmap.conf
share/examples/389-ds-base/dirsrv/config/slapd-collations.conf
share/examples/389-ds-base/dirsrv/schema/00core.ldif
share/examples/389-ds-base/dirsrv/schema/01common.ldif
share/examples/389-ds-base/dirsrv/schema/05rfc2247.ldif
share/examples/389-ds-base/dirsrv/schema/05rfc2927.ldif
share/examples/389-ds-base/dirsrv/schema/10rfc2307.ldif
share/examples/389-ds-base/dirsrv/schema/20subscriber.ldif
share/examples/389-ds-base/dirsrv/schema/25java-object.ldif
share/examples/389-ds-base/dirsrv/schema/28pilot.ldif
share/examples/389-ds-base/dirsrv/schema/30ns-common.ldif
share/examples/389-ds-base/dirsrv/schema/50ns-admin.ldif
share/examples/389-ds-base/dirsrv/schema/50ns-certificate.ldif
share/examples/389-ds-base/dirsrv/schema/50ns-directory.ldif
share/examples/389-ds-base/dirsrv/schema/50ns-mail.ldif
share/examples/389-ds-base/dirsrv/schema/50ns-value.ldif
share/examples/389-ds-base/dirsrv/schema/50ns-web.ldif
share/examples/389-ds-base/dirsrv/schema/60autofs.ldif
share/examples/389-ds-base/dirsrv/schema/60eduperson.ldif
share/examples/389-ds-base/dirsrv/schema/60mozilla.ldif
share/examples/389-ds-base/dirsrv/schema/60nss-ldap.ldif
share/examples/389-ds-base/dirsrv/schema/60pam-plugin.ldif
share/examples/389-ds-base/dirsrv/schema/60pureftpd.ldif
share/examples/389-ds-base/dirsrv/schema/60rfc2739.ldif
share/examples/389-ds-base/dirsrv/schema/60rfc3712.ldif
share/examples/389-ds-base/dirsrv/schema/60sabayon.ldif
share/examples/389-ds-base/dirsrv/schema/60sudo.ldif
share/examples/389-ds-base/dirsrv/schema/60trust.ldif
share/examples/389-ds-base/dirsrv/schema/99user.ldif
share/examples/389-ds-base/rc.d/init.d/dirsrv
