# influxdb/PLIST
@comment $NetBSD$
sbin/influxd
share/examples/rc.d/influxdb
