# bitcoin/PLIST
@comment $NetBSD: PLIST,v 1.4 2014/10/20 09:49:29 noud4 Exp $
bin/bitcoin-qt
sbin/bitcoind
share/applications/bitcoin-qt.desktop
share/examples/bitcoin/bitcoin.conf
share/examples/rc.d/bitcoind
share/pixmaps/bitcoin128.png
