# bitcoin/PLIST
@comment $NetBSD: PLIST,v 1.2 2013/12/21 14:53:09 noud4 Exp $
bin/bitcoin-qt
bin/bitcoind
@comment share/examples/bitcoin/bitcoin.conf
@comment share/examples/rc.d/bitcoind
