# bitcoin/PLIST
@comment $NetBSD: PLIST,v 1.3 2013/12/22 22:25:59 noud4 Exp $
bin/bitcoin-qt
sbin/bitcoind
share/examples/bitcoin/bitcoin.conf
share/examples/rc.d/bitcoind
