# qico/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:03 jeremy-c-reed Exp $
bin/qico
${QCC}bin/qcc
share/examples/qico/qico.conf
share/examples/rc.d/qico
@dirrm share/examples/qico
