# qico/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/10/11 10:45:13 thomasklausner Exp $
${QCC}bin/qcc
bin/qctl
man/man8/qcc.8
man/man8/qctl.8
man/man8/qico.8
sbin/qico
share/examples/qico/qico.conf
share/examples/rc.d/qico
