# qico/PLIST
@comment $NetBSD: PLIST,v 1.4 2007/03/03 15:22:17 shattered Exp $
${QCC}bin/qcc
bin/qctl
man/man8/qcc.8
man/man8/qctl.8
man/man8/qico.8
sbin/qico
share/examples/qico/qico.conf
share/examples/rc.d/qico
@dirrm share/examples/qico
