# jetty/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2006/07/17 20:47:29 unex_linija Exp $
jetty/LICENSE.TXT
jetty/README.TXT
jetty/VERSION.TXT
jetty/etc/LICENSE.javax.servlet.txt
jetty/etc/LICENSE.javax.xml.html
jetty/etc/LICENSE.jsse.txt
jetty/etc/admin.xml
jetty/etc/jetty-jmx.xml
jetty/etc/jetty.xml
jetty/etc/webdefault.xml
jetty/ext/ant.jar
jetty/ext/commons-el.jar
jetty/ext/commons-logging.jar
jetty/ext/jasper-compiler.jar
jetty/ext/jasper-runtime.jar
jetty/ext/mx4j-remote.jar
jetty/ext/mx4j-tools.jar
jetty/ext/mx4j.jar
jetty/ext/xercesImpl.jar
jetty/ext/xml-apis.jar
jetty/ext/xmlParserAPIs.jar
jetty/extra/lib/org.mortbay.jetty-jdk1.2.jar
jetty/extra/lib/org.mortbay.jmx-jdk1.2.jar
jetty/lib/javax.servlet.jar
jetty/lib/org.mortbay.jetty.jar
jetty/lib/org.mortbay.jmx.jar
jetty/start.jar
jetty/stop.jar
jetty/webapps/javadoc.war
jetty/webapps/template/WEB-INF/classes/com/acme/DemoFilter.class
jetty/webapps/template/WEB-INF/classes/com/acme/DumpServlet.class
jetty/webapps/template/WEB-INF/lib/.keepme
jetty/webapps/template/WEB-INF/web.xml
jetty/webapps/template/build.xml
jetty/webapps/template/index.html
jetty/webapps/template/jetty_banner.gif
jetty/webapps/template/src/com/acme/DemoFilter.java
jetty/webapps/template/src/com/acme/DumpServlet.java
share/examples/jetty/admin.xml
share/examples/jetty/jetty-jmx.xml
share/examples/jetty/jetty.xml
share/examples/jetty/webdefault.xml
share/examples/rc.d/jetty
@dirrm share/examples/jetty
@dirrm jetty/webapps/template/src/com/acme
@dirrm jetty/webapps/template/src/com
@dirrm jetty/webapps/template/src
@dirrm jetty/webapps/template/WEB-INF/lib
@dirrm jetty/webapps/template/WEB-INF/classes/com/acme
@dirrm jetty/webapps/template/WEB-INF/classes/com
@dirrm jetty/webapps/template/WEB-INF/classes
@dirrm jetty/webapps/template/WEB-INF
@dirrm jetty/webapps/template
@dirrm jetty/webapps
@exec ${MKDIR} %D/jetty/logs
@dirrm jetty/logs
@dirrm jetty/lib
@dirrm jetty/extra/lib
@dirrm jetty/extra
@dirrm jetty/ext
@dirrm jetty/etc
@dirrm jetty
