# redis/PLIST
@comment $NetBSD: PLIST,v 1.7 2012/02/02 08:52:50 fhajny Exp $
bin/redis-benchmark
bin/redis-check-aof
bin/redis-check-dump
bin/redis-cli
bin/redis-server
share/doc/redis/00-RELEASENOTES
share/doc/redis/BUGS
share/doc/redis/COPYING
share/doc/redis/README
share/examples/rc.d/redis
share/examples/redis/redis.conf.example
share/examples/redis/tests/assets/default.conf
share/examples/redis/tests/helpers/gen_write_load.tcl
share/examples/redis/tests/integration/aof-race.tcl
share/examples/redis/tests/integration/aof.tcl
share/examples/redis/tests/integration/redis-cli.tcl
share/examples/redis/tests/integration/replication-2.tcl
share/examples/redis/tests/integration/replication-3.tcl
share/examples/redis/tests/integration/replication.tcl
share/examples/redis/tests/support/redis.tcl
share/examples/redis/tests/support/server.tcl
share/examples/redis/tests/support/test.tcl
share/examples/redis/tests/support/tmpfile.tcl
share/examples/redis/tests/support/util.tcl
share/examples/redis/tests/test_helper.tcl
share/examples/redis/tests/tmp/.gitignore
share/examples/redis/tests/unit/auth.tcl
share/examples/redis/tests/unit/basic.tcl
share/examples/redis/tests/unit/cas.tcl
share/examples/redis/tests/unit/expire.tcl
share/examples/redis/tests/unit/introspection.tcl
share/examples/redis/tests/unit/maxmemory.tcl
share/examples/redis/tests/unit/other.tcl
share/examples/redis/tests/unit/printver.tcl
share/examples/redis/tests/unit/protocol.tcl
share/examples/redis/tests/unit/pubsub.tcl
share/examples/redis/tests/unit/quit.tcl
share/examples/redis/tests/unit/slowlog.tcl
share/examples/redis/tests/unit/sort.tcl
share/examples/redis/tests/unit/type/hash.tcl
share/examples/redis/tests/unit/type/list-2.tcl
share/examples/redis/tests/unit/type/list-3.tcl
share/examples/redis/tests/unit/type/list-common.tcl
share/examples/redis/tests/unit/type/list.tcl
share/examples/redis/tests/unit/type/set.tcl
share/examples/redis/tests/unit/type/zset.tcl
share/examples/redis/utils/generate-command-help.rb
share/examples/redis/utils/install_server.sh
share/examples/redis/utils/mktarball.sh
share/examples/redis/utils/redis.conf.tpl
share/examples/redis/utils/redis_init_script
share/examples/redis/utils/redis_init_script.tpl
