# nvtv/PLIST
@comment $NetBSD: PLIST,v 1.3 2003/06/07 17:50:19 cubidou Exp $
etc/rc.d/nvtvd
bin/nvtv
bin/nvtvd
share/doc/nvtv/FAQ
share/doc/nvtv/README
share/doc/nvtv/USAGE
share/doc/nvtv/chips.txt
share/doc/nvtv/overview.txt
share/doc/nvtv/timing.txt
man/man1/nvtv.1
man/man8/nvtvd.8
@dirrm share/doc/nvtv
